# tinycorelinux/9.x/x86/tcz/quagga.tcz.list
/usr/local/tce.installed/quagga
/usr/local/sbin/zebra
/usr/local/sbin/watchquagga
/usr/local/sbin/ripngd
/usr/local/sbin/ripd
/usr/local/sbin/pimd
/usr/local/sbin/ospfd
/usr/local/sbin/ospfclient
/usr/local/sbin/ospf6d
/usr/local/sbin/nhrpd
/usr/local/sbin/isisd
/usr/local/sbin/bgpd
/usr/local/lib/libzebra.so.1.0.0
/usr/local/lib/libzebra.so.1
/usr/local/lib/libzebra.so
/usr/local/lib/libquagga_pb.so.0.0.0
/usr/local/lib/libquagga_pb.so.0
/usr/local/lib/libquagga_pb.so
/usr/local/lib/libospfapiclient.so.0.0.0
/usr/local/lib/libospfapiclient.so.0
/usr/local/lib/libospfapiclient.so
/usr/local/lib/libospf.so.0.0.0
/usr/local/lib/libospf.so.0
/usr/local/lib/libospf.so
/usr/local/lib/libfpm_pb.so.0.0.0
/usr/local/lib/libfpm_pb.so.0
/usr/local/lib/libfpm_pb.so
/usr/local/etc/quagga/zebra.conf.sample
/usr/local/etc/quagga/vtysh.conf.sample
/usr/local/etc/quagga/ripngd.conf.sample
/usr/local/etc/quagga/ripd.conf.sample
/usr/local/etc/quagga/pimd.conf.sample
/usr/local/etc/quagga/ospfd.conf.sample
/usr/local/etc/quagga/ospf6d.conf.sample
/usr/local/etc/quagga/isisd.conf.sample
/usr/local/etc/quagga/bgpd.conf.sample2
/usr/local/etc/quagga/bgpd.conf.sample
/usr/local/etc/init.d/quagga
/usr/local/bin/vtysh
/usr/local/bin/test_igmpv3_join
/usr/local/bin/bgp_btoa
