# tinycorelinux/4.x/x86/tcz/aiccu.tcz.list
usr/local/etc/aiccu.conf
usr/local/etc/init.d/aiccu
usr/local/sbin/aiccu
usr/local/share/doc/aiccu/HOWTO
usr/local/share/doc/aiccu/LICENSE
usr/local/share/doc/aiccu/README
