# tinycorelinux/4.x/x86/tcz/anyterm.tcz.list
usr/local/bin/anytermd
usr/local/etc/init.d/anyterm
