# tinycorelinux/4.x/x86/tcz/aprx.tcz.list
usr/local/etc/aprx.conf.sample
usr/local/etc/init.d/aprx
usr/local/sbin/aprx
usr/local/sbin/aprx-stat
usr/local/share/doc/aprx/LICENSE
