# tinycorelinux/4.x/x86/tcz/cherokee.tcz.list
usr/local/bin/CTK-run
usr/local/bin/cget
usr/local/bin/cherokee-admin-launcher
usr/local/bin/cherokee-config
usr/local/bin/cherokee-macos-askpass
usr/local/bin/cherokee-panic
usr/local/bin/cherokee-tweak
usr/local/etc/cherokee/cherokee.conf
usr/local/etc/cherokee/cherokee.conf.perf_sample
usr/local/etc/init.d/cherokee
usr/local/lib/cherokee/libplugin_admin.la
usr/local/lib/cherokee/libplugin_admin.so
usr/local/lib/cherokee/libplugin_and.la
usr/local/lib/cherokee/libplugin_and.so
usr/local/lib/cherokee/libplugin_authlist.la
usr/local/lib/cherokee/libplugin_authlist.so
usr/local/lib/cherokee/libplugin_bind.la
usr/local/lib/cherokee/libplugin_bind.so
usr/local/lib/cherokee/libplugin_cgi.la
usr/local/lib/cherokee/libplugin_cgi.so
usr/local/lib/cherokee/libplugin_combined.la
usr/local/lib/cherokee/libplugin_combined.so
usr/local/lib/cherokee/libplugin_common.la
usr/local/lib/cherokee/libplugin_common.so
usr/local/lib/cherokee/libplugin_custom.la
usr/local/lib/cherokee/libplugin_custom.so
usr/local/lib/cherokee/libplugin_custom_error.la
usr/local/lib/cherokee/libplugin_custom_error.so
usr/local/lib/cherokee/libplugin_deflate.la
usr/local/lib/cherokee/libplugin_deflate.so
usr/local/lib/cherokee/libplugin_directory.la
usr/local/lib/cherokee/libplugin_directory.so
usr/local/lib/cherokee/libplugin_dirlist.la
usr/local/lib/cherokee/libplugin_dirlist.so
usr/local/lib/cherokee/libplugin_drop.la
usr/local/lib/cherokee/libplugin_drop.so
usr/local/lib/cherokee/libplugin_empty_gif.la
usr/local/lib/cherokee/libplugin_empty_gif.so
usr/local/lib/cherokee/libplugin_error_nn.la
usr/local/lib/cherokee/libplugin_error_nn.so
usr/local/lib/cherokee/libplugin_error_redir.la
usr/local/lib/cherokee/libplugin_error_redir.so
usr/local/lib/cherokee/libplugin_evhost.la
usr/local/lib/cherokee/libplugin_evhost.so
usr/local/lib/cherokee/libplugin_exists.la
usr/local/lib/cherokee/libplugin_exists.so
usr/local/lib/cherokee/libplugin_extensions.la
usr/local/lib/cherokee/libplugin_extensions.so
usr/local/lib/cherokee/libplugin_failover.la
usr/local/lib/cherokee/libplugin_failover.so
usr/local/lib/cherokee/libplugin_fcgi.la
usr/local/lib/cherokee/libplugin_fcgi.so
usr/local/lib/cherokee/libplugin_file.la
usr/local/lib/cherokee/libplugin_file.so
usr/local/lib/cherokee/libplugin_from.la
usr/local/lib/cherokee/libplugin_from.so
usr/local/lib/cherokee/libplugin_fullpath.la
usr/local/lib/cherokee/libplugin_fullpath.so
usr/local/lib/cherokee/libplugin_gzip.la
usr/local/lib/cherokee/libplugin_gzip.so
usr/local/lib/cherokee/libplugin_header.la
usr/local/lib/cherokee/libplugin_header.so
usr/local/lib/cherokee/libplugin_htdigest.la
usr/local/lib/cherokee/libplugin_htdigest.so
usr/local/lib/cherokee/libplugin_htpasswd.la
usr/local/lib/cherokee/libplugin_htpasswd.so
usr/local/lib/cherokee/libplugin_ip_hash.la
usr/local/lib/cherokee/libplugin_ip_hash.so
usr/local/lib/cherokee/libplugin_libssl.la
usr/local/lib/cherokee/libplugin_libssl.so
usr/local/lib/cherokee/libplugin_method.la
usr/local/lib/cherokee/libplugin_method.so
usr/local/lib/cherokee/libplugin_ncsa.la
usr/local/lib/cherokee/libplugin_ncsa.so
usr/local/lib/cherokee/libplugin_not.la
usr/local/lib/cherokee/libplugin_not.so
usr/local/lib/cherokee/libplugin_or.la
usr/local/lib/cherokee/libplugin_or.so
usr/local/lib/cherokee/libplugin_plain.la
usr/local/lib/cherokee/libplugin_plain.so
usr/local/lib/cherokee/libplugin_post_report.la
usr/local/lib/cherokee/libplugin_post_report.so
usr/local/lib/cherokee/libplugin_post_track.la
usr/local/lib/cherokee/libplugin_post_track.so
usr/local/lib/cherokee/libplugin_proxy.la
usr/local/lib/cherokee/libplugin_proxy.so
usr/local/lib/cherokee/libplugin_redir.la
usr/local/lib/cherokee/libplugin_redir.so
usr/local/lib/cherokee/libplugin_rehost.la
usr/local/lib/cherokee/libplugin_rehost.so
usr/local/lib/cherokee/libplugin_render_rrd.la
usr/local/lib/cherokee/libplugin_render_rrd.so
usr/local/lib/cherokee/libplugin_request.la
usr/local/lib/cherokee/libplugin_request.so
usr/local/lib/cherokee/libplugin_round_robin.la
usr/local/lib/cherokee/libplugin_round_robin.so
usr/local/lib/cherokee/libplugin_rrd.la
usr/local/lib/cherokee/libplugin_rrd.so
usr/local/lib/cherokee/libplugin_scgi.la
usr/local/lib/cherokee/libplugin_scgi.so
usr/local/lib/cherokee/libplugin_secdownload.la
usr/local/lib/cherokee/libplugin_secdownload.so
usr/local/lib/cherokee/libplugin_server_info.la
usr/local/lib/cherokee/libplugin_server_info.so
usr/local/lib/cherokee/libplugin_ssi.la
usr/local/lib/cherokee/libplugin_ssi.so
usr/local/lib/cherokee/libplugin_streaming.la
usr/local/lib/cherokee/libplugin_streaming.so
usr/local/lib/cherokee/libplugin_target_ip.la
usr/local/lib/cherokee/libplugin_target_ip.so
usr/local/lib/cherokee/libplugin_tls.la
usr/local/lib/cherokee/libplugin_tls.so
usr/local/lib/cherokee/libplugin_url_arg.la
usr/local/lib/cherokee/libplugin_url_arg.so
usr/local/lib/cherokee/libplugin_uwsgi.la
usr/local/lib/cherokee/libplugin_uwsgi.so
usr/local/lib/cherokee/libplugin_v_or.la
usr/local/lib/cherokee/libplugin_v_or.so
usr/local/lib/cherokee/libplugin_wildcard.la
usr/local/lib/cherokee/libplugin_wildcard.so
usr/local/lib/libcherokee-base.la
usr/local/lib/libcherokee-base.so
usr/local/lib/libcherokee-base.so.0
usr/local/lib/libcherokee-base.so.0.0.1
usr/local/lib/libcherokee-client.la
usr/local/lib/libcherokee-client.so
usr/local/lib/libcherokee-client.so.0
usr/local/lib/libcherokee-client.so.0.0.1
usr/local/lib/libcherokee-server.la
usr/local/lib/libcherokee-server.so
usr/local/lib/libcherokee-server.so.0
usr/local/lib/libcherokee-server.so.0.0.1
usr/local/lib/pkgconfig/cherokee.pc
usr/local/sbin/cherokee
usr/local/sbin/cherokee-admin
usr/local/sbin/cherokee-worker
usr/local/sbin/htdigest.py
usr/local/share/cherokee/admin/CTK/CTK/AjaxUpload.py
usr/local/share/cherokee/admin/CTK/CTK/Box.py
usr/local/share/cherokee/admin/CTK/CTK/Button.py
usr/local/share/cherokee/admin/CTK/CTK/Carousel.py
usr/local/share/cherokee/admin/CTK/CTK/Checkbox.py
usr/local/share/cherokee/admin/CTK/CTK/Collapsible.py
usr/local/share/cherokee/admin/CTK/CTK/Combobox.py
usr/local/share/cherokee/admin/CTK/CTK/Config.py
usr/local/share/cherokee/admin/CTK/CTK/Container.py
usr/local/share/cherokee/admin/CTK/CTK/DatePicker.py
usr/local/share/cherokee/admin/CTK/CTK/Dialog.py
usr/local/share/cherokee/admin/CTK/CTK/Downloader.py
usr/local/share/cherokee/admin/CTK/CTK/Druid.py
usr/local/share/cherokee/admin/CTK/CTK/HTTP.py
usr/local/share/cherokee/admin/CTK/CTK/Help.py
usr/local/share/cherokee/admin/CTK/CTK/HiddenField.py
usr/local/share/cherokee/admin/CTK/CTK/Image.py
usr/local/share/cherokee/admin/CTK/CTK/Indenter.py
usr/local/share/cherokee/admin/CTK/CTK/Init.py
usr/local/share/cherokee/admin/CTK/CTK/JS.py
usr/local/share/cherokee/admin/CTK/CTK/Link.py
usr/local/share/cherokee/admin/CTK/CTK/List.py
usr/local/share/cherokee/admin/CTK/CTK/MailHTML.py
usr/local/share/cherokee/admin/CTK/CTK/Notice.py
usr/local/share/cherokee/admin/CTK/CTK/Page.py
usr/local/share/cherokee/admin/CTK/CTK/PageCleaner.py
usr/local/share/cherokee/admin/CTK/CTK/Paginator.py
usr/local/share/cherokee/admin/CTK/CTK/Plugin.py
usr/local/share/cherokee/admin/CTK/CTK/Post.py
usr/local/share/cherokee/admin/CTK/CTK/ProgressBar.py
usr/local/share/cherokee/admin/CTK/CTK/PropsTable.py
usr/local/share/cherokee/admin/CTK/CTK/Proxy.py
usr/local/share/cherokee/admin/CTK/CTK/Radio.py
usr/local/share/cherokee/admin/CTK/CTK/RawHTML.py
usr/local/share/cherokee/admin/CTK/CTK/Refreshable.py
usr/local/share/cherokee/admin/CTK/CTK/Server.py
usr/local/share/cherokee/admin/CTK/CTK/SortableList.py
usr/local/share/cherokee/admin/CTK/CTK/StarRating.py
usr/local/share/cherokee/admin/CTK/CTK/Submitter.py
usr/local/share/cherokee/admin/CTK/CTK/Tab.py
usr/local/share/cherokee/admin/CTK/CTK/Table.py
usr/local/share/cherokee/admin/CTK/CTK/Template.py
usr/local/share/cherokee/admin/CTK/CTK/TextArea.py
usr/local/share/cherokee/admin/CTK/CTK/TextField.py
usr/local/share/cherokee/admin/CTK/CTK/ToggleButton.py
usr/local/share/cherokee/admin/CTK/CTK/Uploader.py
usr/local/share/cherokee/admin/CTK/CTK/Widget.py
usr/local/share/cherokee/admin/CTK/CTK/XMLRPCProxy.py
usr/local/share/cherokee/admin/CTK/CTK/__init__.py
usr/local/share/cherokee/admin/CTK/CTK/consts.py
usr/local/share/cherokee/admin/CTK/CTK/i18n.py
usr/local/share/cherokee/admin/CTK/CTK/iPhoneToggle.py
usr/local/share/cherokee/admin/CTK/CTK/json_embedded.py
usr/local/share/cherokee/admin/CTK/CTK/pyscgi.py
usr/local/share/cherokee/admin/CTK/CTK/util.py
usr/local/share/cherokee/admin/CTK/static/css/CTK.css
usr/local/share/cherokee/admin/CTK/static/css/datepicker.css
usr/local/share/cherokee/admin/CTK/static/css/jquery-ui-1.7.2.custom.css
usr/local/share/cherokee/admin/CTK/static/css/jquery.ibutton.css
usr/local/share/cherokee/admin/CTK/static/images/arrow_ns.png
usr/local/share/cherokee/admin/CTK/static/images/bg-dialog-buttonpane.png
usr/local/share/cherokee/admin/CTK/static/images/bg-dialog-titlebar-close.gif
usr/local/share/cherokee/admin/CTK/static/images/bg-dialog-titlebar.png
usr/local/share/cherokee/admin/CTK/static/images/bg-dialog.png
usr/local/share/cherokee/admin/CTK/static/images/bg-filter.png
usr/local/share/cherokee/admin/CTK/static/images/bg-sel-actions-sel.png
usr/local/share/cherokee/admin/CTK/static/images/bg-sel-actions.png
usr/local/share/cherokee/admin/CTK/static/images/carousel-left-arrow.png
usr/local/share/cherokee/admin/CTK/static/images/carousel-right-arrow.png
usr/local/share/cherokee/admin/CTK/static/images/del-hover.png
usr/local/share/cherokee/admin/CTK/static/images/del.png
usr/local/share/cherokee/admin/CTK/static/images/dialog-error.png
usr/local/share/cherokee/admin/CTK/static/images/dialog-information.png
usr/local/share/cherokee/admin/CTK/static/images/dialog-warning.png
usr/local/share/cherokee/admin/CTK/static/images/forbid.png
usr/local/share/cherokee/admin/CTK/static/images/ibutton-slider-default.png
usr/local/share/cherokee/admin/CTK/static/images/input-bg.png
usr/local/share/cherokee/admin/CTK/static/images/jquery.ui.stars.gif
usr/local/share/cherokee/admin/CTK/static/images/loader.gif
usr/local/share/cherokee/admin/CTK/static/images/loading.gif
usr/local/share/cherokee/admin/CTK/static/images/off.png
usr/local/share/cherokee/admin/CTK/static/images/offline.png
usr/local/share/cherokee/admin/CTK/static/images/on.png
usr/local/share/cherokee/admin/CTK/static/images/online.png
usr/local/share/cherokee/admin/CTK/static/images/tabbg.png
usr/local/share/cherokee/admin/CTK/static/images/tick.png
usr/local/share/cherokee/admin/CTK/static/images/ui-bg_gloss-wave_35_f6a828_500x100.png
usr/local/share/cherokee/admin/CTK/static/images/ui-bg_highlight-soft_100_eeeeee_1x100.png
usr/local/share/cherokee/admin/CTK/static/images/ui-icons_222222_256x240.png
usr/local/share/cherokee/admin/CTK/static/images/ui-icons_2e83ff_256x240.png
usr/local/share/cherokee/admin/CTK/static/images/ui-icons_454545_256x240.png
usr/local/share/cherokee/admin/CTK/static/images/ui-icons_888888_256x240.png
usr/local/share/cherokee/admin/CTK/static/images/ui-icons_cd0a0a_256x240.png
usr/local/share/cherokee/admin/CTK/static/images/uploadify.cancel.png
usr/local/share/cherokee/admin/CTK/static/js/Carousel.js
usr/local/share/cherokee/admin/CTK/static/js/Help.js
usr/local/share/cherokee/admin/CTK/static/js/StarRating.js
usr/local/share/cherokee/admin/CTK/static/js/Submitter.js
usr/local/share/cherokee/admin/CTK/static/js/ajaxupload.3.6.js
usr/local/share/cherokee/admin/CTK/static/js/common.js
usr/local/share/cherokee/admin/CTK/static/js/jquery-1.3.2.js
usr/local/share/cherokee/admin/CTK/static/js/jquery-1.3.2.min.js
usr/local/share/cherokee/admin/CTK/static/js/jquery-ui-1.7.2.custom.min.js
usr/local/share/cherokee/admin/CTK/static/js/jquery-ui-1.7.2.js
usr/local/share/cherokee/admin/CTK/static/js/jquery.cookie.js
usr/local/share/cherokee/admin/CTK/static/js/jquery.form-defaults.js
usr/local/share/cherokee/admin/CTK/static/js/jquery.ibutton.js
usr/local/share/cherokee/admin/CTK/static/js/jquery.tablednd_0_5.js
usr/local/share/cherokee/admin/CTK/static/js/jquery.uploadProgress.js
usr/local/share/cherokee/admin/cherokee.conf.sample
usr/local/share/cherokee/admin/icons/cherokee-admin-launcher-128.png
usr/local/share/cherokee/admin/icons/cherokee-admin-launcher-16.png
usr/local/share/cherokee/admin/icons/cherokee-admin-launcher-256.png
usr/local/share/cherokee/admin/icons/cherokee-admin-launcher-32.png
usr/local/share/cherokee/admin/icons/cherokee-admin-launcher-48.png
usr/local/share/cherokee/admin/icons/cherokee-admin-launcher-512.png
usr/local/share/cherokee/admin/icons/cherokee-admin-launcher-96.png
usr/local/share/cherokee/admin/icons/cherokee-admin-launcher.svg
usr/local/share/cherokee/admin/market/CommandProgress.py
usr/local/share/cherokee/admin/market/Distro.py
usr/local/share/cherokee/admin/market/Install.py
usr/local/share/cherokee/admin/market/InstallUtil.py
usr/local/share/cherokee/admin/market/Install_Log.py
usr/local/share/cherokee/admin/market/Library.py
usr/local/share/cherokee/admin/market/Maintenance.py
usr/local/share/cherokee/admin/market/Menu.py
usr/local/share/cherokee/admin/market/PageApp.py
usr/local/share/cherokee/admin/market/PageCategory.py
usr/local/share/cherokee/admin/market/PageIndex.py
usr/local/share/cherokee/admin/market/PageSearch.py
usr/local/share/cherokee/admin/market/Report.py
usr/local/share/cherokee/admin/market/Review.py
usr/local/share/cherokee/admin/market/Util.py
usr/local/share/cherokee/admin/market/__init__.py
usr/local/share/cherokee/admin/market/ows_consts.py
usr/local/share/cherokee/admin/performance.conf.sample
usr/local/share/cherokee/admin/plugins/admin.py
usr/local/share/cherokee/admin/plugins/authlist.py
usr/local/share/cherokee/admin/plugins/bind.py
usr/local/share/cherokee/admin/plugins/cgi.py
usr/local/share/cherokee/admin/plugins/common.py
usr/local/share/cherokee/admin/plugins/custom_error.py
usr/local/share/cherokee/admin/plugins/dbslayer.py
usr/local/share/cherokee/admin/plugins/deflate.py
usr/local/share/cherokee/admin/plugins/directory.py
usr/local/share/cherokee/admin/plugins/dirlist.py
usr/local/share/cherokee/admin/plugins/drop.py
usr/local/share/cherokee/admin/plugins/empty_gif.py
usr/local/share/cherokee/admin/plugins/error_nn.py
usr/local/share/cherokee/admin/plugins/error_redir.py
usr/local/share/cherokee/admin/plugins/evhost.py
usr/local/share/cherokee/admin/plugins/exists.py
usr/local/share/cherokee/admin/plugins/extensions.py
usr/local/share/cherokee/admin/plugins/failover.py
usr/local/share/cherokee/admin/plugins/fcgi.py
usr/local/share/cherokee/admin/plugins/file.py
usr/local/share/cherokee/admin/plugins/from.py
usr/local/share/cherokee/admin/plugins/fullpath.py
usr/local/share/cherokee/admin/plugins/geoip.py
usr/local/share/cherokee/admin/plugins/gzip.py
usr/local/share/cherokee/admin/plugins/header.py
usr/local/share/cherokee/admin/plugins/htdigest.py
usr/local/share/cherokee/admin/plugins/htpasswd.py
usr/local/share/cherokee/admin/plugins/ip_hash.py
usr/local/share/cherokee/admin/plugins/ldap.py
usr/local/share/cherokee/admin/plugins/method.py
usr/local/share/cherokee/admin/plugins/mysql.py
usr/local/share/cherokee/admin/plugins/pam.py
usr/local/share/cherokee/admin/plugins/plain.py
usr/local/share/cherokee/admin/plugins/post_report.py
usr/local/share/cherokee/admin/plugins/post_track.py
usr/local/share/cherokee/admin/plugins/proxy.py
usr/local/share/cherokee/admin/plugins/redir.py
usr/local/share/cherokee/admin/plugins/rehost.py
usr/local/share/cherokee/admin/plugins/request.py
usr/local/share/cherokee/admin/plugins/round_robin.py
usr/local/share/cherokee/admin/plugins/rrd.py
usr/local/share/cherokee/admin/plugins/scgi.py
usr/local/share/cherokee/admin/plugins/secdownload.py
usr/local/share/cherokee/admin/plugins/server_info.py
usr/local/share/cherokee/admin/plugins/ssi.py
usr/local/share/cherokee/admin/plugins/streaming.py
usr/local/share/cherokee/admin/plugins/target_ip.py
usr/local/share/cherokee/admin/plugins/tls.py
usr/local/share/cherokee/admin/plugins/url_arg.py
usr/local/share/cherokee/admin/plugins/uwsgi.py
usr/local/share/cherokee/admin/plugins/wildcard.py
usr/local/share/cherokee/admin/wizards/List.py
usr/local/share/cherokee/admin/wizards/alfresco.py
usr/local/share/cherokee/admin/wizards/coldfusion.py
usr/local/share/cherokee/admin/wizards/concrete5.py
usr/local/share/cherokee/admin/wizards/django.py
usr/local/share/cherokee/admin/wizards/flcache.py
usr/local/share/cherokee/admin/wizards/glassfish.py
usr/local/share/cherokee/admin/wizards/hotlinking.py
usr/local/share/cherokee/admin/wizards/icons.py
usr/local/share/cherokee/admin/wizards/mailman.py
usr/local/share/cherokee/admin/wizards/mono.py
usr/local/share/cherokee/admin/wizards/php.py
usr/local/share/cherokee/admin/wizards/rails.py
usr/local/share/cherokee/admin/wizards/redirect.py
usr/local/share/cherokee/admin/wizards/rtorrent.py
usr/local/share/cherokee/admin/wizards/ssl_test.py
usr/local/share/cherokee/admin/wizards/static.py
usr/local/share/cherokee/admin/wizards/streaming.py
usr/local/share/cherokee/admin/wizards/symfony.py
usr/local/share/cherokee/admin/wizards/uwsgi.py
usr/local/share/cherokee/admin/wizards/zend.py
usr/local/share/cherokee/deps/cgi.deps
usr/local/share/cherokee/deps/combined.deps
usr/local/share/cherokee/deps/common.deps
usr/local/share/cherokee/deps/error_redir.deps
usr/local/share/cherokee/deps/fcgi.deps
usr/local/share/cherokee/deps/post_report.deps
usr/local/share/cherokee/deps/render_rrd.deps
usr/local/share/cherokee/deps/scgi.deps
usr/local/share/cherokee/deps/secdownload.deps
usr/local/share/cherokee/deps/streaming.deps
usr/local/share/cherokee/deps/uwsgi.deps
usr/local/share/cherokee/icons/README.txt
usr/local/share/cherokee/icons/arrow_turn_left.png
usr/local/share/cherokee/icons/blank.png
usr/local/share/cherokee/icons/bomb.png
usr/local/share/cherokee/icons/camera.png
usr/local/share/cherokee/icons/cd.png
usr/local/share/cherokee/icons/color_wheel.png
usr/local/share/cherokee/icons/control_play.png
usr/local/share/cherokee/icons/css.png
usr/local/share/cherokee/icons/cup.png
usr/local/share/cherokee/icons/email.png
usr/local/share/cherokee/icons/film.png
usr/local/share/cherokee/icons/folder.png
usr/local/share/cherokee/icons/font.png
usr/local/share/cherokee/icons/html.png
usr/local/share/cherokee/icons/music.png
usr/local/share/cherokee/icons/package.png
usr/local/share/cherokee/icons/page_white.png
usr/local/share/cherokee/icons/page_white_acrobat.png
usr/local/share/cherokee/icons/page_white_actionscript.png
usr/local/share/cherokee/icons/page_white_c.png
usr/local/share/cherokee/icons/page_white_go.png
usr/local/share/cherokee/icons/page_white_office.png
usr/local/share/cherokee/icons/page_white_php.png
usr/local/share/cherokee/icons/page_white_text.png
usr/local/share/cherokee/icons/printer.png
usr/local/share/cherokee/icons/ruby.png
usr/local/share/cherokee/icons/script.png
usr/local/share/cherokee/themes/default/entry.html
usr/local/share/cherokee/themes/default/footer.html
usr/local/share/cherokee/themes/default/header.html
usr/local/share/cherokee/themes/default/logo.png
usr/local/share/cherokee/themes/default/theme.css
usr/local/share/cherokee/themes/firefox3/entry.html
usr/local/share/cherokee/themes/firefox3/footer.html
usr/local/share/cherokee/themes/firefox3/header.html
usr/local/share/cherokee/themes/firefox3/theme.css
usr/local/share/cherokee/themes/plain/entry.html
usr/local/share/cherokee/themes/plain/footer.html
usr/local/share/cherokee/themes/plain/header.html
usr/local/share/cherokee/themes/plain/theme.css
usr/local/share/cherokee/themes/white/entry.html
usr/local/share/cherokee/themes/white/footer.html
usr/local/share/cherokee/themes/white/header.html
usr/local/share/cherokee/themes/white/theme.css
var/www/images/cherokee-logo.png
var/www/images/default-bg.png
var/www/images/favicon.ico
var/www/images/powered_by_cherokee.png
var/www/index.html
