# tinycorelinux/4.x/x86/tcz/chrony.tcz.list
usr/local/bin/chronyc
usr/local/etc/chrony/chrony.conf.sample
usr/local/etc/chrony/chrony.keys.sample
usr/local/etc/init.d/chrony
usr/local/sbin/chronyd
usr/local/tce.installed/chrony
