# tinycorelinux/4.x/x86/tcz/cups-browsed.tcz.list
usr/local/etc/cups/templates/cups-browsed.conf
usr/local/etc/init.d/2b
usr/local/etc/init.d/browsed
usr/local/sbin/cups-browsed
usr/local/share/doc/cups-browsed/AUTHORS
usr/local/share/doc/cups-browsed/COPYING
usr/local/share/doc/cups-browsed/README
usr/local/tce.installed/cups-browsed
