# tinycorelinux/4.x/x86/tcz/dansguardian.tcz.list
usr/local/etc/dansguardian/authplugins/ident.conf
usr/local/etc/dansguardian/authplugins/ip.conf
usr/local/etc/dansguardian/authplugins/proxy-basic.conf
usr/local/etc/dansguardian/authplugins/proxy-digest.conf
usr/local/etc/dansguardian/downloadmanagers/default.conf
usr/local/etc/dansguardian/downloadmanagers/fancy.conf
usr/local/etc/dansguardian/lists/authplugins/ipgroups
usr/local/etc/dansguardian/lists/bannedextensionlist
usr/local/etc/dansguardian/lists/bannediplist
usr/local/etc/dansguardian/lists/bannedmimetypelist
usr/local/etc/dansguardian/lists/bannedphraselist
usr/local/etc/dansguardian/lists/bannedregexpheaderlist
usr/local/etc/dansguardian/lists/bannedregexpurllist
usr/local/etc/dansguardian/lists/bannedsitelist
usr/local/etc/dansguardian/lists/bannedurllist
usr/local/etc/dansguardian/lists/blacklists/ads/domains
usr/local/etc/dansguardian/lists/blacklists/ads/urls
usr/local/etc/dansguardian/lists/contentregexplist
usr/local/etc/dansguardian/lists/downloadmanagers/managedextensionlist
usr/local/etc/dansguardian/lists/downloadmanagers/managedmimetypelist
usr/local/etc/dansguardian/lists/exceptionextensionlist
usr/local/etc/dansguardian/lists/exceptionfilesitelist
usr/local/etc/dansguardian/lists/exceptionfileurllist
usr/local/etc/dansguardian/lists/exceptioniplist
usr/local/etc/dansguardian/lists/exceptionmimetypelist
usr/local/etc/dansguardian/lists/exceptionphraselist
usr/local/etc/dansguardian/lists/exceptionregexpurllist
usr/local/etc/dansguardian/lists/exceptionsitelist
usr/local/etc/dansguardian/lists/exceptionurllist
usr/local/etc/dansguardian/lists/filtergroupslist
usr/local/etc/dansguardian/lists/greysitelist
usr/local/etc/dansguardian/lists/greyurllist
usr/local/etc/dansguardian/lists/headerregexplist
usr/local/etc/dansguardian/lists/logregexpurllist
usr/local/etc/dansguardian/lists/logsitelist
usr/local/etc/dansguardian/lists/logurllist
usr/local/etc/dansguardian/lists/phraselists/badwords/weighted_dutch
usr/local/etc/dansguardian/lists/phraselists/badwords/weighted_french
usr/local/etc/dansguardian/lists/phraselists/badwords/weighted_german
usr/local/etc/dansguardian/lists/phraselists/badwords/weighted_portuguese
usr/local/etc/dansguardian/lists/phraselists/badwords/weighted_spanish
usr/local/etc/dansguardian/lists/phraselists/chat/weighted
usr/local/etc/dansguardian/lists/phraselists/chat/weighted_italian
usr/local/etc/dansguardian/lists/phraselists/conspiracy/weighted
usr/local/etc/dansguardian/lists/phraselists/domainsforsale/weighted
usr/local/etc/dansguardian/lists/phraselists/drugadvocacy/weighted
usr/local/etc/dansguardian/lists/phraselists/forums/weighted
usr/local/etc/dansguardian/lists/phraselists/gambling/banned
usr/local/etc/dansguardian/lists/phraselists/gambling/banned_portuguese
usr/local/etc/dansguardian/lists/phraselists/gambling/weighted
usr/local/etc/dansguardian/lists/phraselists/gambling/weighted_portuguese
usr/local/etc/dansguardian/lists/phraselists/games/weighted
usr/local/etc/dansguardian/lists/phraselists/goodphrases/exception
usr/local/etc/dansguardian/lists/phraselists/goodphrases/exception_email
usr/local/etc/dansguardian/lists/phraselists/goodphrases/weighted_general
usr/local/etc/dansguardian/lists/phraselists/goodphrases/weighted_general_danish
usr/local/etc/dansguardian/lists/phraselists/goodphrases/weighted_general_dutch
usr/local/etc/dansguardian/lists/phraselists/goodphrases/weighted_general_malay
usr/local/etc/dansguardian/lists/phraselists/goodphrases/weighted_general_polish
usr/local/etc/dansguardian/lists/phraselists/goodphrases/weighted_general_portuguese
usr/local/etc/dansguardian/lists/phraselists/goodphrases/weighted_general_swedish
usr/local/etc/dansguardian/lists/phraselists/goodphrases/weighted_news
usr/local/etc/dansguardian/lists/phraselists/googlesearches/banned
usr/local/etc/dansguardian/lists/phraselists/gore/weighted
usr/local/etc/dansguardian/lists/phraselists/gore/weighted_portuguese
usr/local/etc/dansguardian/lists/phraselists/idtheft/weighted
usr/local/etc/dansguardian/lists/phraselists/illegaldrugs/banned
usr/local/etc/dansguardian/lists/phraselists/illegaldrugs/weighted
usr/local/etc/dansguardian/lists/phraselists/illegaldrugs/weighted_portuguese
usr/local/etc/dansguardian/lists/phraselists/intolerance/banned_portuguese
usr/local/etc/dansguardian/lists/phraselists/intolerance/weighted
usr/local/etc/dansguardian/lists/phraselists/intolerance/weighted_portuguese
usr/local/etc/dansguardian/lists/phraselists/legaldrugs/weighted
usr/local/etc/dansguardian/lists/phraselists/malware/weighted
usr/local/etc/dansguardian/lists/phraselists/music/weighted
usr/local/etc/dansguardian/lists/phraselists/news/weighted
usr/local/etc/dansguardian/lists/phraselists/nudism/weighted
usr/local/etc/dansguardian/lists/phraselists/peer2peer/weighted
usr/local/etc/dansguardian/lists/phraselists/personals/weighted
usr/local/etc/dansguardian/lists/phraselists/personals/weighted_portuguese
usr/local/etc/dansguardian/lists/phraselists/pornography/banned
usr/local/etc/dansguardian/lists/phraselists/pornography/banned_portuguese
usr/local/etc/dansguardian/lists/phraselists/pornography/weighted
usr/local/etc/dansguardian/lists/phraselists/pornography/weighted_chinese
usr/local/etc/dansguardian/lists/phraselists/pornography/weighted_danish
usr/local/etc/dansguardian/lists/phraselists/pornography/weighted_dutch
usr/local/etc/dansguardian/lists/phraselists/pornography/weighted_french
usr/local/etc/dansguardian/lists/phraselists/pornography/weighted_german
usr/local/etc/dansguardian/lists/phraselists/pornography/weighted_italian
usr/local/etc/dansguardian/lists/phraselists/pornography/weighted_japanese
usr/local/etc/dansguardian/lists/phraselists/pornography/weighted_malay
usr/local/etc/dansguardian/lists/phraselists/pornography/weighted_norwegian
usr/local/etc/dansguardian/lists/phraselists/pornography/weighted_polish
usr/local/etc/dansguardian/lists/phraselists/pornography/weighted_portuguese
usr/local/etc/dansguardian/lists/phraselists/pornography/weighted_russian
usr/local/etc/dansguardian/lists/phraselists/pornography/weighted_spanish
usr/local/etc/dansguardian/lists/phraselists/pornography/weighted_swedish
usr/local/etc/dansguardian/lists/phraselists/proxies/weighted
usr/local/etc/dansguardian/lists/phraselists/rta/banned
usr/local/etc/dansguardian/lists/phraselists/safelabel/banned
usr/local/etc/dansguardian/lists/phraselists/secretsocieties/weighted
usr/local/etc/dansguardian/lists/phraselists/sport/weighted
usr/local/etc/dansguardian/lists/phraselists/translation/weighted
usr/local/etc/dansguardian/lists/phraselists/travel/weighted
usr/local/etc/dansguardian/lists/phraselists/upstreamfilter/weighted
usr/local/etc/dansguardian/lists/phraselists/violence/weighted
usr/local/etc/dansguardian/lists/phraselists/violence/weighted_portuguese
usr/local/etc/dansguardian/lists/phraselists/warezhacking/weighted
usr/local/etc/dansguardian/lists/phraselists/weapons/weighted
usr/local/etc/dansguardian/lists/phraselists/weapons/weighted_portuguese
usr/local/etc/dansguardian/lists/phraselists/webmail/weighted
usr/local/etc/dansguardian/lists/pics
usr/local/etc/dansguardian/lists/urlregexplist
usr/local/etc/dansguardian/lists/weightedphraselist
usr/local/etc/init.d/dansguardian
usr/local/sbin/dansguardian
usr/local/share/dansguardian/dansguardian.pl
usr/local/share/dansguardian/files/dansguardian.conf
usr/local/share/dansguardian/files/dansguardianf1.conf
usr/local/share/dansguardian/languages/arspanish/fancydmtemplate.html
usr/local/share/dansguardian/languages/arspanish/messages
usr/local/share/dansguardian/languages/arspanish/template.html
usr/local/share/dansguardian/languages/bulgarian/fancydmtemplate.html
usr/local/share/dansguardian/languages/bulgarian/messages
usr/local/share/dansguardian/languages/bulgarian/template.html
usr/local/share/dansguardian/languages/chinesebig5/fancydmtemplate.html
usr/local/share/dansguardian/languages/chinesebig5/messages
usr/local/share/dansguardian/languages/chinesebig5/template.html
usr/local/share/dansguardian/languages/chinesegb2312/fancydmtemplate.html
usr/local/share/dansguardian/languages/chinesegb2312/messages
usr/local/share/dansguardian/languages/chinesegb2312/template.html
usr/local/share/dansguardian/languages/czech/fancydmtemplate.html
usr/local/share/dansguardian/languages/czech/messages
usr/local/share/dansguardian/languages/czech/template.html
usr/local/share/dansguardian/languages/danish/fancydmtemplate.html
usr/local/share/dansguardian/languages/danish/messages
usr/local/share/dansguardian/languages/danish/template.html
usr/local/share/dansguardian/languages/dutch/fancydmtemplate.html
usr/local/share/dansguardian/languages/dutch/messages
usr/local/share/dansguardian/languages/dutch/template.html
usr/local/share/dansguardian/languages/french/fancydmtemplate.html
usr/local/share/dansguardian/languages/french/messages
usr/local/share/dansguardian/languages/french/template.html
usr/local/share/dansguardian/languages/german/fancydmtemplate.html
usr/local/share/dansguardian/languages/german/messages
usr/local/share/dansguardian/languages/german/template.html
usr/local/share/dansguardian/languages/hebrew/fancydmtemplate.html
usr/local/share/dansguardian/languages/hebrew/messages
usr/local/share/dansguardian/languages/hebrew/template.html
usr/local/share/dansguardian/languages/hungarian/fancydmtemplate.html
usr/local/share/dansguardian/languages/hungarian/messages
usr/local/share/dansguardian/languages/hungarian/template.html
usr/local/share/dansguardian/languages/indonesian/fancydmtemplate.html
usr/local/share/dansguardian/languages/indonesian/messages
usr/local/share/dansguardian/languages/indonesian/template.html
usr/local/share/dansguardian/languages/italian/fancydmtemplate.html
usr/local/share/dansguardian/languages/italian/messages
usr/local/share/dansguardian/languages/italian/template.html
usr/local/share/dansguardian/languages/japanese/fancydmtemplate.html
usr/local/share/dansguardian/languages/japanese/messages
usr/local/share/dansguardian/languages/japanese/template.html
usr/local/share/dansguardian/languages/lithuanian/fancydmtemplate.html
usr/local/share/dansguardian/languages/lithuanian/messages
usr/local/share/dansguardian/languages/lithuanian/template.html
usr/local/share/dansguardian/languages/malay/fancydmtemplate.html
usr/local/share/dansguardian/languages/malay/messages
usr/local/share/dansguardian/languages/malay/template.html
usr/local/share/dansguardian/languages/mxspanish/fancydmtemplate.html
usr/local/share/dansguardian/languages/mxspanish/messages
usr/local/share/dansguardian/languages/mxspanish/template.html
usr/local/share/dansguardian/languages/polish/fancydmtemplate.html
usr/local/share/dansguardian/languages/polish/messages
usr/local/share/dansguardian/languages/polish/template.html
usr/local/share/dansguardian/languages/portuguese/fancydmtemplate.html
usr/local/share/dansguardian/languages/portuguese/messages
usr/local/share/dansguardian/languages/portuguese/template.html
usr/local/share/dansguardian/languages/ptbrazilian/fancydmtemplate.html
usr/local/share/dansguardian/languages/ptbrazilian/messages
usr/local/share/dansguardian/languages/ptbrazilian/template.html
usr/local/share/dansguardian/languages/russian-1251/fancydmtemplate.html
usr/local/share/dansguardian/languages/russian-1251/messages
usr/local/share/dansguardian/languages/russian-1251/template.html
usr/local/share/dansguardian/languages/russian-koi8-r/fancydmtemplate.html
usr/local/share/dansguardian/languages/russian-koi8-r/messages
usr/local/share/dansguardian/languages/russian-koi8-r/template.html
usr/local/share/dansguardian/languages/slovak/fancydmtemplate.html
usr/local/share/dansguardian/languages/slovak/messages
usr/local/share/dansguardian/languages/slovak/template.html
usr/local/share/dansguardian/languages/spanish/fancydmtemplate.html
usr/local/share/dansguardian/languages/spanish/messages
usr/local/share/dansguardian/languages/spanish/template.html
usr/local/share/dansguardian/languages/swedish/fancydmtemplate.html
usr/local/share/dansguardian/languages/swedish/messages
usr/local/share/dansguardian/languages/swedish/template.html
usr/local/share/dansguardian/languages/turkish/fancydmtemplate.html
usr/local/share/dansguardian/languages/turkish/messages
usr/local/share/dansguardian/languages/turkish/template.html
usr/local/share/dansguardian/languages/ukenglish/fancydmtemplate.html
usr/local/share/dansguardian/languages/ukenglish/messages
usr/local/share/dansguardian/languages/ukenglish/template.html
usr/local/share/dansguardian/scripts/bsd-init
usr/local/share/dansguardian/scripts/dansguardian
usr/local/share/dansguardian/scripts/logrotation
usr/local/share/dansguardian/scripts/solaris-init
usr/local/share/dansguardian/scripts/systemv-init
usr/local/share/dansguardian/transparent1x1.gif
usr/local/tce.installed/dansguardian
