# tinycorelinux/4.x/x86/tcz/drizzle-server.tcz.list
usr/local/etc/init.d/drizzle-server
usr/local/lib/drizzle7/libarchive_plugin.so
usr/local/lib/drizzle7/libascii_plugin.so
usr/local/lib/drizzle7/libauth_file_plugin.so
usr/local/lib/drizzle7/libauth_http_plugin.so
usr/local/lib/drizzle7/libauth_ldap_plugin.so
usr/local/lib/drizzle7/libauth_pam_plugin.so
usr/local/lib/drizzle7/libauth_test_plugin.so
usr/local/lib/drizzle7/libbenchmark_plugin.so
usr/local/lib/drizzle7/libblackhole_plugin.so
usr/local/lib/drizzle7/libcharlength_plugin.so
usr/local/lib/drizzle7/libcompression_plugin.so
usr/local/lib/drizzle7/libconnection_id_plugin.so
usr/local/lib/drizzle7/libcrc32_plugin.so
usr/local/lib/drizzle7/libdebug_plugin.so
usr/local/lib/drizzle7/libdefault_replicator_plugin.so
usr/local/lib/drizzle7/liberrmsg_stderr_plugin.so
usr/local/lib/drizzle7/libfiltered_replicator_plugin.so
usr/local/lib/drizzle7/libfunction_dictionary_plugin.so
usr/local/lib/drizzle7/libhello_events_plugin.so
usr/local/lib/drizzle7/libhello_world_plugin.so
usr/local/lib/drizzle7/libhex_functions_plugin.so
usr/local/lib/drizzle7/liblength_plugin.so
usr/local/lib/drizzle7/liblogging_query_plugin.so
usr/local/lib/drizzle7/liblogging_stats_plugin.so
usr/local/lib/drizzle7/libmd5_plugin.so
usr/local/lib/drizzle7/libmulti_thread_plugin.so
usr/local/lib/drizzle7/libpbms_plugin.so
usr/local/lib/drizzle7/libperformance_dictionary_plugin.so
usr/local/lib/drizzle7/librand_function_plugin.so
usr/local/lib/drizzle7/libreverse_function_plugin.so
usr/local/lib/drizzle7/librot13_plugin.so
usr/local/lib/drizzle7/libshow_schema_proto_plugin.so
usr/local/lib/drizzle7/libshutdown_function_plugin.so
usr/local/lib/drizzle7/libsimple_user_policy_plugin.so
usr/local/lib/drizzle7/libsleep_plugin.so
usr/local/lib/drizzle7/libstorage_engine_api_tester_plugin.so
usr/local/lib/drizzle7/libsubstr_functions_plugin.so
usr/local/lib/drizzle7/libsyslog_plugin.so
usr/local/lib/drizzle7/libtableprototester_plugin.so
usr/local/lib/drizzle7/libtrigger_dictionary_plugin.so
usr/local/lib/drizzle7/libutility_dictionary_plugin.so
usr/local/lib/drizzle7/libuuid_function_plugin.so
usr/local/lib/drizzle7/libversion_plugin.so
usr/local/sbin/drizzled
usr/local/sbin/drizzled7
usr/local/share/man/man8/drizzled.8
usr/local/share/man/man8/drizzled7.8
usr/local/tce.installed/drizzle-server
