# tinycorelinux/4.x/x86/tcz/ebtables.tcz.list
usr/local/etc/ethertypes
usr/local/etc/rc.d/init.d/ebtables
usr/local/etc/sysconfig/ebtables-config
usr/local/lib/libebt_802_3.so
usr/local/lib/libebt_among.so
usr/local/lib/libebt_arp.so
usr/local/lib/libebt_arpreply.so
usr/local/lib/libebt_ip.so
usr/local/lib/libebt_ip6.so
usr/local/lib/libebt_limit.so
usr/local/lib/libebt_log.so
usr/local/lib/libebt_mark.so
usr/local/lib/libebt_mark_m.so
usr/local/lib/libebt_nat.so
usr/local/lib/libebt_nflog.so
usr/local/lib/libebt_pkttype.so
usr/local/lib/libebt_redirect.so
usr/local/lib/libebt_standard.so
usr/local/lib/libebt_stp.so
usr/local/lib/libebt_ulog.so
usr/local/lib/libebt_vlan.so
usr/local/lib/libebtable_broute.so
usr/local/lib/libebtable_filter.so
usr/local/lib/libebtable_nat.so
usr/local/lib/libebtc.so
usr/local/sbin/ebtables
usr/local/sbin/ebtables-restore
usr/local/sbin/ebtables-save
