# tinycorelinux/4.x/x86/tcz/edna.tcz.list
usr/local/bin/MP3Info.py
usr/local/bin/edna.py
usr/local/bin/ezt.py
usr/local/bin/id3reader.py
usr/local/etc/init.d/edna
usr/local/share/applications/edna.desktop
usr/local/share/applications/edna~1.desktop
usr/local/share/applications/edna~2.desktop
usr/local/share/edna/conf/edna.conf
usr/local/share/edna/resources/folder.png
usr/local/share/edna/resources/play.png
usr/local/share/edna/resources/save.png
usr/local/share/edna/resources/text.png
usr/local/share/edna/templates/default.ezt
usr/local/share/edna/templates/default_complex.ezt
usr/local/share/edna/templates/stats.ezt
usr/local/share/edna/templates/style-xml.ezt
usr/local/tce.installed/edna
