# tinycorelinux/4.x/x86/tcz/firebird-2.5-classic.tcz.list
etc/profile.d/firebird-tools.sh
usr/local/bin/fb_inet_server
usr/local/etc/firebird-2.5-classic/UDF
usr/local/etc/firebird-2.5-classic/bin
usr/local/etc/firebird-2.5-classic/firebird.log
usr/local/etc/firebird-2.5-classic/firebird.msg
usr/local/etc/firebird-2.5-classic/intl/fbintl
usr/local/etc/firebird-2.5-classic/intl/fbintl.conf
usr/local/etc/firebird-2.5-classic/lib
usr/local/etc/firebird-2.5-classic/plugins
usr/local/etc/firebird-2.5-classic/security2.fdb
usr/local/etc/init.d/firebird-2.5-classic
usr/local/lib/firebird-2.5/UDF/fbudf.so
usr/local/lib/firebird-2.5/UDF/ib_udf.so
usr/local/lib/firebird-2.5/aliases.conf.default
usr/local/lib/firebird-2.5/fbtrace.conf.default
usr/local/lib/firebird-2.5/firebird.conf.default
usr/local/lib/firebird-2.5/firebird.msg
usr/local/lib/firebird-2.5/intl/fbintl
usr/local/lib/firebird-2.5/intl/fbintl.conf
usr/local/lib/firebird-2.5/plugins/libfbtrace.so
usr/local/lib/firebird-2.5/security2.fdb.default
usr/local/lib/libfbembed.so
usr/local/lib/libfbembed.so.2.5
usr/local/lib/libfbembed.so.2.5.0
usr/local/lib/libib_util.so
usr/local/tce.installed/firebird-2.5-classic
