# tinycorelinux/4.x/x86/tcz/firebird-2.5-super.tcz.list
etc/profile.d/firebird-tools.sh
usr/local/bin/fbguard
usr/local/bin/fbserver
usr/local/etc/firebird-2.5-super/UDF
usr/local/etc/firebird-2.5-super/bin
usr/local/etc/firebird-2.5-super/firebird.log
usr/local/etc/firebird-2.5-super/firebird.msg
usr/local/etc/firebird-2.5-super/intl/fbintl
usr/local/etc/firebird-2.5-super/intl/fbintl.conf
usr/local/etc/firebird-2.5-super/lib
usr/local/etc/firebird-2.5-super/plugins
usr/local/etc/firebird-2.5-super/security2.fdb
usr/local/etc/init.d/firebird-2.5-super
usr/local/lib/firebird-2.5/UDF/fbudf.so
usr/local/lib/firebird-2.5/UDF/ib_udf.so
usr/local/lib/firebird-2.5/aliases.conf.default
usr/local/lib/firebird-2.5/fbtrace.conf.default
usr/local/lib/firebird-2.5/firebird.conf.default
usr/local/lib/firebird-2.5/firebird.msg
usr/local/lib/firebird-2.5/intl/fbintl
usr/local/lib/firebird-2.5/intl/fbintl.conf
usr/local/lib/firebird-2.5/plugins/libfbtrace.so
usr/local/lib/firebird-2.5/security2.fdb.default
usr/local/lib/libib_util.so
usr/local/tce.installed/firebird-2.5-super
