# tinycorelinux/4.x/x86/tcz/freeradius.tcz.list
usr/local/bin/rad_counter
usr/local/bin/radclient
usr/local/bin/radconf2xml
usr/local/bin/radcrypt
usr/local/bin/radeapclient
usr/local/bin/radlast
usr/local/bin/radsniff
usr/local/bin/radsqlrelay
usr/local/bin/radtest
usr/local/bin/radwho
usr/local/bin/radzap
usr/local/bin/rlm_ippool_tool
usr/local/bin/smbencrypt
usr/local/etc/init.d/freeradius
usr/local/etc/raddb-skel.tgz
usr/local/lib/rlm_acct_unique-2.2.0.so
usr/local/lib/rlm_acct_unique.so
usr/local/lib/rlm_acctlog-2.2.0.so
usr/local/lib/rlm_acctlog.so
usr/local/lib/rlm_always-2.2.0.so
usr/local/lib/rlm_always.so
usr/local/lib/rlm_attr_filter-2.2.0.so
usr/local/lib/rlm_attr_filter.so
usr/local/lib/rlm_attr_rewrite-2.2.0.so
usr/local/lib/rlm_attr_rewrite.so
usr/local/lib/rlm_cache-2.2.0.so
usr/local/lib/rlm_cache.so
usr/local/lib/rlm_caching-2.2.0.so
usr/local/lib/rlm_caching.so
usr/local/lib/rlm_chap-2.2.0.so
usr/local/lib/rlm_chap.so
usr/local/lib/rlm_checkval-2.2.0.so
usr/local/lib/rlm_checkval.so
usr/local/lib/rlm_copy_packet-2.2.0.so
usr/local/lib/rlm_copy_packet.so
usr/local/lib/rlm_counter-2.2.0.so
usr/local/lib/rlm_counter.so
usr/local/lib/rlm_cram-2.2.0.so
usr/local/lib/rlm_cram.so
usr/local/lib/rlm_detail-2.2.0.so
usr/local/lib/rlm_detail.so
usr/local/lib/rlm_digest-2.2.0.so
usr/local/lib/rlm_digest.so
usr/local/lib/rlm_dynamic_clients-2.2.0.so
usr/local/lib/rlm_dynamic_clients.so
usr/local/lib/rlm_eap-2.2.0.so
usr/local/lib/rlm_eap.so
usr/local/lib/rlm_eap_gtc-2.2.0.so
usr/local/lib/rlm_eap_gtc.so
usr/local/lib/rlm_eap_leap-2.2.0.so
usr/local/lib/rlm_eap_leap.so
usr/local/lib/rlm_eap_md5-2.2.0.so
usr/local/lib/rlm_eap_md5.so
usr/local/lib/rlm_eap_mschapv2-2.2.0.so
usr/local/lib/rlm_eap_mschapv2.so
usr/local/lib/rlm_eap_peap-2.2.0.so
usr/local/lib/rlm_eap_peap.so
usr/local/lib/rlm_eap_sim-2.2.0.so
usr/local/lib/rlm_eap_sim.so
usr/local/lib/rlm_eap_tls-2.2.0.so
usr/local/lib/rlm_eap_tls.so
usr/local/lib/rlm_eap_ttls-2.2.0.so
usr/local/lib/rlm_eap_ttls.so
usr/local/lib/rlm_example-2.2.0.so
usr/local/lib/rlm_example.so
usr/local/lib/rlm_exec-2.2.0.so
usr/local/lib/rlm_exec.so
usr/local/lib/rlm_expiration-2.2.0.so
usr/local/lib/rlm_expiration.so
usr/local/lib/rlm_expr-2.2.0.so
usr/local/lib/rlm_expr.so
usr/local/lib/rlm_fastusers-2.2.0.so
usr/local/lib/rlm_fastusers.so
usr/local/lib/rlm_files-2.2.0.so
usr/local/lib/rlm_files.so
usr/local/lib/rlm_ippool-2.2.0.so
usr/local/lib/rlm_ippool.so
usr/local/lib/rlm_jradius-2.2.0.so
usr/local/lib/rlm_jradius.so
usr/local/lib/rlm_krb5-2.2.0.so
usr/local/lib/rlm_krb5.so
usr/local/lib/rlm_ldap-2.2.0.so
usr/local/lib/rlm_ldap.so
usr/local/lib/rlm_linelog-2.2.0.so
usr/local/lib/rlm_linelog.so
usr/local/lib/rlm_logintime-2.2.0.so
usr/local/lib/rlm_logintime.so
usr/local/lib/rlm_mschap-2.2.0.so
usr/local/lib/rlm_mschap.so
usr/local/lib/rlm_otp-2.2.0.so
usr/local/lib/rlm_otp.so
usr/local/lib/rlm_pam-2.2.0.so
usr/local/lib/rlm_pam.so
usr/local/lib/rlm_pap-2.2.0.so
usr/local/lib/rlm_pap.so
usr/local/lib/rlm_passwd-2.2.0.so
usr/local/lib/rlm_passwd.so
usr/local/lib/rlm_perl-2.2.0.so
usr/local/lib/rlm_perl.so
usr/local/lib/rlm_policy-2.2.0.so
usr/local/lib/rlm_policy.so
usr/local/lib/rlm_preprocess-2.2.0.so
usr/local/lib/rlm_preprocess.so
usr/local/lib/rlm_protocol_filter-2.2.0.so
usr/local/lib/rlm_protocol_filter.so
usr/local/lib/rlm_radutmp-2.2.0.so
usr/local/lib/rlm_radutmp.so
usr/local/lib/rlm_realm-2.2.0.so
usr/local/lib/rlm_realm.so
usr/local/lib/rlm_replicate-2.2.0.so
usr/local/lib/rlm_replicate.so
usr/local/lib/rlm_sim_files-2.2.0.so
usr/local/lib/rlm_sim_files.so
usr/local/lib/rlm_smsotp-2.2.0.so
usr/local/lib/rlm_smsotp.so
usr/local/lib/rlm_soh-2.2.0.so
usr/local/lib/rlm_soh.so
usr/local/lib/rlm_unix-2.2.0.so
usr/local/lib/rlm_unix.so
usr/local/lib/rlm_wimax-2.2.0.so
usr/local/lib/rlm_wimax.so
usr/local/sbin/checkrad
usr/local/sbin/raddebug
usr/local/sbin/radiusd
usr/local/sbin/radmin
usr/local/sbin/radwatch
usr/local/sbin/rc.radiusd
usr/local/share/doc/License/freeradius-copyright.txt
usr/local/share/doc/License/freeradius.txt
usr/local/share/freeradius/dictionary
usr/local/share/freeradius/dictionary.3com
usr/local/share/freeradius/dictionary.3gpp
usr/local/share/freeradius/dictionary.3gpp2
usr/local/share/freeradius/dictionary.acc
usr/local/share/freeradius/dictionary.acme
usr/local/share/freeradius/dictionary.airespace
usr/local/share/freeradius/dictionary.alcatel
usr/local/share/freeradius/dictionary.alcatel.esam
usr/local/share/freeradius/dictionary.alcatel.sr
usr/local/share/freeradius/dictionary.alteon
usr/local/share/freeradius/dictionary.altiga
usr/local/share/freeradius/dictionary.alvarion
usr/local/share/freeradius/dictionary.apc
usr/local/share/freeradius/dictionary.aptis
usr/local/share/freeradius/dictionary.aruba
usr/local/share/freeradius/dictionary.ascend
usr/local/share/freeradius/dictionary.ascend.illegal
usr/local/share/freeradius/dictionary.asn
usr/local/share/freeradius/dictionary.audiocodes
usr/local/share/freeradius/dictionary.avaya
usr/local/share/freeradius/dictionary.azaire
usr/local/share/freeradius/dictionary.bay
usr/local/share/freeradius/dictionary.bintec
usr/local/share/freeradius/dictionary.bristol
usr/local/share/freeradius/dictionary.bt
usr/local/share/freeradius/dictionary.cablelabs
usr/local/share/freeradius/dictionary.cabletron
usr/local/share/freeradius/dictionary.camiant
usr/local/share/freeradius/dictionary.chillispot
usr/local/share/freeradius/dictionary.cisco
usr/local/share/freeradius/dictionary.cisco.bbsm
usr/local/share/freeradius/dictionary.cisco.vpn3000
usr/local/share/freeradius/dictionary.cisco.vpn5000
usr/local/share/freeradius/dictionary.clavister
usr/local/share/freeradius/dictionary.colubris
usr/local/share/freeradius/dictionary.columbia_university
usr/local/share/freeradius/dictionary.compat
usr/local/share/freeradius/dictionary.cosine
usr/local/share/freeradius/dictionary.dhcp
usr/local/share/freeradius/dictionary.digium
usr/local/share/freeradius/dictionary.eltex
usr/local/share/freeradius/dictionary.epygi
usr/local/share/freeradius/dictionary.ericsson
usr/local/share/freeradius/dictionary.erx
usr/local/share/freeradius/dictionary.extreme
usr/local/share/freeradius/dictionary.f5
usr/local/share/freeradius/dictionary.fortinet
usr/local/share/freeradius/dictionary.foundry
usr/local/share/freeradius/dictionary.freeradius
usr/local/share/freeradius/dictionary.freeradius.internal
usr/local/share/freeradius/dictionary.freeswitch
usr/local/share/freeradius/dictionary.gandalf
usr/local/share/freeradius/dictionary.garderos
usr/local/share/freeradius/dictionary.gemtek
usr/local/share/freeradius/dictionary.h3c
usr/local/share/freeradius/dictionary.hp
usr/local/share/freeradius/dictionary.huawei
usr/local/share/freeradius/dictionary.iea
usr/local/share/freeradius/dictionary.infonet
usr/local/share/freeradius/dictionary.ipunplugged
usr/local/share/freeradius/dictionary.issanni
usr/local/share/freeradius/dictionary.itk
usr/local/share/freeradius/dictionary.jradius
usr/local/share/freeradius/dictionary.juniper
usr/local/share/freeradius/dictionary.karlnet
usr/local/share/freeradius/dictionary.lancom
usr/local/share/freeradius/dictionary.livingston
usr/local/share/freeradius/dictionary.localweb
usr/local/share/freeradius/dictionary.lucent
usr/local/share/freeradius/dictionary.manzara
usr/local/share/freeradius/dictionary.merit
usr/local/share/freeradius/dictionary.microsoft
usr/local/share/freeradius/dictionary.mikrotik
usr/local/share/freeradius/dictionary.motorola
usr/local/share/freeradius/dictionary.motorola.illegal
usr/local/share/freeradius/dictionary.motorola.wimax
usr/local/share/freeradius/dictionary.navini
usr/local/share/freeradius/dictionary.netscreen
usr/local/share/freeradius/dictionary.networkphysics
usr/local/share/freeradius/dictionary.nexans
usr/local/share/freeradius/dictionary.nokia
usr/local/share/freeradius/dictionary.nokia.conflict
usr/local/share/freeradius/dictionary.nomadix
usr/local/share/freeradius/dictionary.nortel
usr/local/share/freeradius/dictionary.ntua
usr/local/share/freeradius/dictionary.openser
usr/local/share/freeradius/dictionary.packeteer
usr/local/share/freeradius/dictionary.paloalto
usr/local/share/freeradius/dictionary.patton
usr/local/share/freeradius/dictionary.propel
usr/local/share/freeradius/dictionary.prosoft
usr/local/share/freeradius/dictionary.purewave
usr/local/share/freeradius/dictionary.quiconnect
usr/local/share/freeradius/dictionary.quintum
usr/local/share/freeradius/dictionary.redback
usr/local/share/freeradius/dictionary.redcreek
usr/local/share/freeradius/dictionary.rfc2865
usr/local/share/freeradius/dictionary.rfc2866
usr/local/share/freeradius/dictionary.rfc2867
usr/local/share/freeradius/dictionary.rfc2868
usr/local/share/freeradius/dictionary.rfc2869
usr/local/share/freeradius/dictionary.rfc3162
usr/local/share/freeradius/dictionary.rfc3576
usr/local/share/freeradius/dictionary.rfc3580
usr/local/share/freeradius/dictionary.rfc4072
usr/local/share/freeradius/dictionary.rfc4372
usr/local/share/freeradius/dictionary.rfc4603
usr/local/share/freeradius/dictionary.rfc4675
usr/local/share/freeradius/dictionary.rfc4679
usr/local/share/freeradius/dictionary.rfc4818
usr/local/share/freeradius/dictionary.rfc4849
usr/local/share/freeradius/dictionary.rfc5090
usr/local/share/freeradius/dictionary.rfc5176
usr/local/share/freeradius/dictionary.rfc5447
usr/local/share/freeradius/dictionary.rfc5580
usr/local/share/freeradius/dictionary.rfc5607
usr/local/share/freeradius/dictionary.rfc5904
usr/local/share/freeradius/dictionary.riverstone
usr/local/share/freeradius/dictionary.roaringpenguin
usr/local/share/freeradius/dictionary.shasta
usr/local/share/freeradius/dictionary.shiva
usr/local/share/freeradius/dictionary.siemens
usr/local/share/freeradius/dictionary.slipstream
usr/local/share/freeradius/dictionary.sofaware
usr/local/share/freeradius/dictionary.sonicwall
usr/local/share/freeradius/dictionary.springtide
usr/local/share/freeradius/dictionary.starent
usr/local/share/freeradius/dictionary.starent.vsa1
usr/local/share/freeradius/dictionary.symbol
usr/local/share/freeradius/dictionary.t_systems_nova
usr/local/share/freeradius/dictionary.telebit
usr/local/share/freeradius/dictionary.telkom
usr/local/share/freeradius/dictionary.terena
usr/local/share/freeradius/dictionary.trapeze
usr/local/share/freeradius/dictionary.tropos
usr/local/share/freeradius/dictionary.ukerna
usr/local/share/freeradius/dictionary.unix
usr/local/share/freeradius/dictionary.usr
usr/local/share/freeradius/dictionary.usr.illegal
usr/local/share/freeradius/dictionary.utstarcom
usr/local/share/freeradius/dictionary.valemount
usr/local/share/freeradius/dictionary.versanet
usr/local/share/freeradius/dictionary.vqp
usr/local/share/freeradius/dictionary.walabi
usr/local/share/freeradius/dictionary.waverider
usr/local/share/freeradius/dictionary.wichorus
usr/local/share/freeradius/dictionary.wimax
usr/local/share/freeradius/dictionary.wimax.wichorus
usr/local/share/freeradius/dictionary.wispr
usr/local/share/freeradius/dictionary.xedia
usr/local/share/freeradius/dictionary.xylan
usr/local/share/freeradius/dictionary.zyxel
usr/local/tce.installed/freeradius
