# tinycorelinux/4.x/x86/tcz/gpsd.tcz.list
lib/udev/gpsd.hotplug
lib/udev/rules.d/99-gpsd-usb.rules
usr/local/etc/gpsd.rc
usr/local/etc/init.d/gpsd
usr/local/lib/libgps.so
usr/local/lib/libgps.so.20
usr/local/lib/libgps.so.20.0
usr/local/lib/libgps.so.20.0.0
usr/local/lib/libgpsd.so
usr/local/lib/libgpsd.so.20
usr/local/lib/libgpsd.so.20.0
usr/local/lib/libgpsd.so.20.0.0
usr/local/sbin/gpsd
usr/local/sbin/gpsdctl
usr/local/share/man/man8/gpsd.8
usr/local/share/man/man8/gpsdctl.8
