# tinycorelinux/4.x/x86/tcz/hal.tcz.list
usr/local/bin/hal-device
usr/local/bin/hal-disable-polling
usr/local/bin/hal-find-by-capability
usr/local/bin/hal-find-by-property
usr/local/bin/hal-get-property
usr/local/bin/hal-is-caller-locked-out
usr/local/bin/hal-lock
usr/local/bin/hal-set-property
usr/local/bin/hal-setup-keymap
usr/local/bin/lshal
usr/local/etc/dbus-1/system.d/hal.conf
usr/local/etc/init.d/haldaemon
usr/local/lib/libhal-storage.so
usr/local/lib/libhal-storage.so.1
usr/local/lib/libhal-storage.so.1.0.0
usr/local/lib/libhal.so
usr/local/lib/libhal.so.1
usr/local/lib/libhal.so.1.0.0
usr/local/lib/udev/rules.d/90-hal.rules
usr/local/libexec/hal-storage-cleanup-all-mountpoints
usr/local/libexec/hal-storage-cleanup-mountpoint
usr/local/libexec/hal-storage-closetray
usr/local/libexec/hal-storage-eject
usr/local/libexec/hal-storage-mount
usr/local/libexec/hal-storage-unmount
usr/local/libexec/hal-system-power-pm-is-supported
usr/local/libexec/hal-system-power-pmu
usr/local/libexec/hal-system-setserial
usr/local/libexec/hal-system-sonypic
usr/local/libexec/hald-addon-acpi
usr/local/libexec/hald-addon-cpufreq
usr/local/libexec/hald-addon-generic-backlight
usr/local/libexec/hald-addon-hid-ups
usr/local/libexec/hald-addon-imac-backlight
usr/local/libexec/hald-addon-input
usr/local/libexec/hald-addon-ipw-killswitch
usr/local/libexec/hald-addon-leds
usr/local/libexec/hald-addon-macbook-backlight
usr/local/libexec/hald-addon-macbookpro-backlight
usr/local/libexec/hald-addon-pmu
usr/local/libexec/hald-addon-rfkill-killswitch
usr/local/libexec/hald-addon-storage
usr/local/libexec/hald-generate-fdi-cache
usr/local/libexec/hald-probe-hiddev
usr/local/libexec/hald-probe-ieee1394-unit
usr/local/libexec/hald-probe-input
usr/local/libexec/hald-probe-net-bluetooth
usr/local/libexec/hald-probe-pc-floppy
usr/local/libexec/hald-probe-printer
usr/local/libexec/hald-probe-serial
usr/local/libexec/hald-probe-smbios
usr/local/libexec/hald-probe-storage
usr/local/libexec/hald-probe-video4linux
usr/local/libexec/hald-probe-volume
usr/local/libexec/hald-runner
usr/local/libexec/scripts/hal-dockstation-undock
usr/local/libexec/scripts/hal-functions
usr/local/libexec/scripts/hal-luks-remove
usr/local/libexec/scripts/hal-luks-setup
usr/local/libexec/scripts/hal-luks-teardown
usr/local/libexec/scripts/hal-system-killswitch-get-power
usr/local/libexec/scripts/hal-system-killswitch-set-power
usr/local/libexec/scripts/hal-system-lcd-get-brightness
usr/local/libexec/scripts/hal-system-lcd-set-brightness
usr/local/libexec/scripts/hal-system-power-hibernate
usr/local/libexec/scripts/hal-system-power-reboot
usr/local/libexec/scripts/hal-system-power-set-power-save
usr/local/libexec/scripts/hal-system-power-shutdown
usr/local/libexec/scripts/hal-system-power-suspend
usr/local/libexec/scripts/hal-system-power-suspend-hybrid
usr/local/libexec/scripts/hal-system-wol-enable
usr/local/libexec/scripts/hal-system-wol-enabled
usr/local/libexec/scripts/hal-system-wol-supported
usr/local/libexec/scripts/linux/hal-dockstation-undock-linux
usr/local/libexec/scripts/linux/hal-luks-remove-linux
usr/local/libexec/scripts/linux/hal-luks-setup-linux
usr/local/libexec/scripts/linux/hal-luks-teardown-linux
usr/local/libexec/scripts/linux/hal-system-killswitch-get-power-linux
usr/local/libexec/scripts/linux/hal-system-killswitch-set-power-linux
usr/local/libexec/scripts/linux/hal-system-lcd-get-brightness-linux
usr/local/libexec/scripts/linux/hal-system-lcd-set-brightness-linux
usr/local/libexec/scripts/linux/hal-system-power-hibernate-linux
usr/local/libexec/scripts/linux/hal-system-power-reboot-linux
usr/local/libexec/scripts/linux/hal-system-power-set-power-save-linux
usr/local/libexec/scripts/linux/hal-system-power-shutdown-linux
usr/local/libexec/scripts/linux/hal-system-power-suspend-hybrid-linux
usr/local/libexec/scripts/linux/hal-system-power-suspend-linux
usr/local/libexec/scripts/linux/hal-system-wol-enable-linux
usr/local/libexec/scripts/linux/hal-system-wol-enabled-linux
usr/local/libexec/scripts/linux/hal-system-wol-linux
usr/local/libexec/scripts/linux/hal-system-wol-supported-linux
usr/local/sbin/hald
usr/local/share/hal/fdi/fdi.dtd
usr/local/share/hal/fdi/policy/10osvendor/10-cpufreq.fdi
usr/local/share/hal/fdi/policy/10osvendor/10-dockstation.fdi
usr/local/share/hal/fdi/policy/10osvendor/10-imac-backlight.fdi
usr/local/share/hal/fdi/policy/10osvendor/10-input-policy.fdi
usr/local/share/hal/fdi/policy/10osvendor/10-keymap.fdi
usr/local/share/hal/fdi/policy/10osvendor/10-laptop-panel-mgmt-policy.fdi
usr/local/share/hal/fdi/policy/10osvendor/10-leds.fdi
usr/local/share/hal/fdi/policy/10osvendor/10-macbook-backlight.fdi
usr/local/share/hal/fdi/policy/10osvendor/10-macbookpro-utils.fdi
usr/local/share/hal/fdi/policy/10osvendor/10-power-mgmt-policy.fdi
usr/local/share/hal/fdi/policy/10osvendor/10-rfkill-switch.fdi
usr/local/share/hal/fdi/policy/10osvendor/10-tabletPCs.fdi
usr/local/share/hal/fdi/policy/10osvendor/10-x11-input.fdi
usr/local/share/hal/fdi/policy/10osvendor/15-storage-luks.fdi
usr/local/share/hal/fdi/policy/10osvendor/20-storage-methods.fdi
usr/local/share/hal/fdi/policy/10osvendor/30-wol.fdi
