# tinycorelinux/4.x/x86/tcz/isc-dhcp.tcz.list
usr/local/bin/omshell
usr/local/etc/init.d/isc-dhcpd
usr/local/sbin/dhclient
usr/local/sbin/dhclient-script
usr/local/sbin/dhcpd
usr/local/sbin/dhcrelay
