# tinycorelinux/4.x/x86/tcz/laptop-mode-tools.tcz.list
usr/local/etc/init.d/laptop-mode
usr/local/sbin/laptop_mode
usr/local/sbin/lm-profiler
usr/local/sbin/lm-syslog-setup
usr/local/share/laptop-mode-tools/files/acpi/actions/lm_ac_adapter.sh
usr/local/share/laptop-mode-tools/files/acpi/actions/lm_battery.sh
usr/local/share/laptop-mode-tools/files/acpi/actions/lm_lid.sh
usr/local/share/laptop-mode-tools/files/acpi/events/lm_ac_adapter
usr/local/share/laptop-mode-tools/files/acpi/events/lm_battery
usr/local/share/laptop-mode-tools/files/acpi/events/lm_lid
usr/local/share/laptop-mode-tools/files/laptop-mode/conf.d/ac97-powersave.conf
usr/local/share/laptop-mode-tools/files/laptop-mode/conf.d/auto-hibernate.conf
usr/local/share/laptop-mode-tools/files/laptop-mode/conf.d/battery-level-polling.conf
usr/local/share/laptop-mode-tools/files/laptop-mode/conf.d/bluetooth.conf
usr/local/share/laptop-mode-tools/files/laptop-mode/conf.d/configuration-file-control.conf
usr/local/share/laptop-mode-tools/files/laptop-mode/conf.d/cpufreq.conf
usr/local/share/laptop-mode-tools/files/laptop-mode/conf.d/dpms-standby.conf
usr/local/share/laptop-mode-tools/files/laptop-mode/conf.d/eee-superhe.conf
usr/local/share/laptop-mode-tools/files/laptop-mode/conf.d/ethernet.conf
usr/local/share/laptop-mode-tools/files/laptop-mode/conf.d/hal-polling.conf
usr/local/share/laptop-mode-tools/files/laptop-mode/conf.d/intel-hda-powersave.conf
usr/local/share/laptop-mode-tools/files/laptop-mode/conf.d/intel-sata-powermgmt.conf
usr/local/share/laptop-mode-tools/files/laptop-mode/conf.d/lcd-brightness.conf
usr/local/share/laptop-mode-tools/files/laptop-mode/conf.d/sched-mc-power-savings.conf
usr/local/share/laptop-mode-tools/files/laptop-mode/conf.d/sched-smt-power-savings.conf
usr/local/share/laptop-mode-tools/files/laptop-mode/conf.d/start-stop-programs.conf
usr/local/share/laptop-mode-tools/files/laptop-mode/conf.d/terminal-blanking.conf
usr/local/share/laptop-mode-tools/files/laptop-mode/conf.d/usb-autosuspend.conf
usr/local/share/laptop-mode-tools/files/laptop-mode/conf.d/video-out.conf
usr/local/share/laptop-mode-tools/files/laptop-mode/conf.d/wireless-ipw-power.conf
usr/local/share/laptop-mode-tools/files/laptop-mode/conf.d/wireless-iwl-power.conf
usr/local/share/laptop-mode-tools/files/laptop-mode/laptop-mode.conf
usr/local/share/laptop-mode-tools/files/laptop-mode/lm-profiler.conf
usr/local/share/laptop-mode-tools/module-helpers/lm-polling-daemon
usr/local/share/laptop-mode-tools/modules/ac97-powersave
usr/local/share/laptop-mode-tools/modules/battery-level-polling
usr/local/share/laptop-mode-tools/modules/bluetooth
usr/local/share/laptop-mode-tools/modules/configuration-file-control
usr/local/share/laptop-mode-tools/modules/cpufreq
usr/local/share/laptop-mode-tools/modules/dpms-standby
usr/local/share/laptop-mode-tools/modules/eee-superhe
usr/local/share/laptop-mode-tools/modules/ethernet
usr/local/share/laptop-mode-tools/modules/hal-polling
usr/local/share/laptop-mode-tools/modules/hdparm
usr/local/share/laptop-mode-tools/modules/intel-hda-powersave
usr/local/share/laptop-mode-tools/modules/intel-sata-powermgmt
usr/local/share/laptop-mode-tools/modules/laptop-mode
usr/local/share/laptop-mode-tools/modules/lcd-brightness
usr/local/share/laptop-mode-tools/modules/sched-mc-power-savings
usr/local/share/laptop-mode-tools/modules/sched-smt-power-savings
usr/local/share/laptop-mode-tools/modules/start-stop-programs
usr/local/share/laptop-mode-tools/modules/syslog-conf
usr/local/share/laptop-mode-tools/modules/terminal-blanking
usr/local/share/laptop-mode-tools/modules/usb-autosuspend
usr/local/share/laptop-mode-tools/modules/video-out
usr/local/share/laptop-mode-tools/modules/wireless-ipw-power
usr/local/share/laptop-mode-tools/modules/wireless-iwl-power
usr/local/tce.installed/laptop-mode-tools
