# tinycorelinux/4.x/x86/tcz/lcdproc-server.tcz.list
usr/local/etc/LCDd.conf.sample
usr/local/etc/init.d/LCDd
usr/local/lib/lcdproc/curses.so
usr/local/lib/lcdproc/hd44780.so
usr/local/sbin/LCDd
usr/local/tce.installed/lcdproc-server
