# tinycorelinux/4.x/x86/tcz/monkey.tcz.list
usr/local/bin/monkey
usr/local/etc/init.d/banana
usr/local/share/applications/monkey.desktop
usr/local/share/applications/monkey~1.desktop
usr/local/share/applications/monkey~2.desktop
usr/local/share/monkey/cgi-bin/test.pl
usr/local/share/monkey/conf/modules.conf
usr/local/share/monkey/conf/monkey.conf
usr/local/share/monkey/conf/monkey.deny
usr/local/share/monkey/conf/monkey.mime
usr/local/share/monkey/htdocs/htdocs.tar.gz
usr/local/tce.installed/monkey
