# tinycorelinux/4.x/x86/tcz/NetworkManager.tcz.list
etc/udev/rules.d/77-nm-olpc-mesh.rules
usr/local/bin/nm-online
usr/local/bin/nm-tool
usr/local/bin/nmcli
usr/local/etc/NetworkManager/nm-system-settings.conf
usr/local/etc/dbus-1/system.d/NetworkManager.conf
usr/local/etc/dbus-1/system.d/nm-avahi-autoipd.conf
usr/local/etc/dbus-1/system.d/nm-dhcp-client.conf
usr/local/etc/dbus-1/system.d/nm-dispatcher.conf
usr/local/etc/default/NetworkManager
usr/local/etc/init.d/NetworkManager
usr/local/lib/NetworkManager/libnm-settings-plugin-keyfile.so
usr/local/lib/libnm-glib-vpn.so
usr/local/lib/libnm-glib-vpn.so.1
usr/local/lib/libnm-glib-vpn.so.1.0.0
usr/local/lib/libnm-glib.so
usr/local/lib/libnm-glib.so.2
usr/local/lib/libnm-glib.so.2.4.3
usr/local/lib/libnm-util.so
usr/local/lib/libnm-util.so.1
usr/local/lib/libnm-util.so.1.6.1
usr/local/lib/pppd/2.4.4/nm-pppd-plugin.so
usr/local/libexec/nm-avahi-autoipd.action
usr/local/libexec/nm-crash-logger
usr/local/libexec/nm-dhcp-client.action
usr/local/libexec/nm-dispatcher.action
usr/local/sbin/NetworkManager
usr/local/share/NetworkManager/gdb-cmd
usr/local/share/dbus-1/system-services/org.freedesktop.nm_dispatcher.service
usr/local/share/polkit-1/actions/org.freedesktop.NetworkManager.policy
usr/local/share/polkit-1/actions/org.freedesktop.network-manager-settings.system.policy
usr/local/tce.installed/NetworkManager
