# tinycorelinux/4.x/x86/tcz/ntp.tcz.list
usr/local/bin/ntp-keygen
usr/local/bin/ntp-wait
usr/local/bin/ntpd
usr/local/bin/ntpdate
usr/local/bin/ntpdc
usr/local/bin/ntpq
usr/local/bin/ntptime
usr/local/bin/ntptrace
usr/local/bin/sntp
usr/local/bin/tickadj
usr/local/etc/init.d/ntpd
usr/local/share/man/man1/ntp-keygen.1
usr/local/share/man/man1/ntpd.1
usr/local/share/man/man1/ntpdc.1
usr/local/share/man/man1/ntpq.1
usr/local/share/man/man1/ntpsnmpd.1
usr/local/share/man/man1/sntp.1
usr/local/tce.installed/ntp
