# tinycorelinux/4.x/x86/tcz/openslp.tcz.list
usr/local/bin/slptool
usr/local/etc/init.d/slpd
usr/local/etc/slp.conf
usr/local/etc/slp.reg
usr/local/etc/slp.spi
usr/local/lib/libslp.so
usr/local/lib/libslp.so.1
usr/local/lib/libslp.so.1.0.1
usr/local/sbin/slpd
