# tinycorelinux/4.x/x86/tcz/open-vm-tools-3.0.21-tinycore.tcz.list
etc/init.d/open-vm-tools
etc/release
etc/vmware-tools/locations
etc/vmware-tools/plugins
etc/vmware-tools/poweroff-vm-default
etc/vmware-tools/poweron-vm-default
etc/vmware-tools/resume-vm-default
etc/vmware-tools/scripts/vmware/network
etc/vmware-tools/statechange.subr
etc/vmware-tools/suspend-vm-default
etc/vmware-tools/vm-support
usr/local/bin/vmtoolsd
usr/local/bin/vmware-checkvm
usr/local/bin/vmware-hgfsclient
usr/local/bin/vmware-rpctool
usr/local/bin/vmware-toolbox-cmd
usr/local/bin/vmware-user
usr/local/bin/vmware-user-suid-wrapper
usr/local/bin/vmware-vmblock-fuse
usr/local/bin/vmware-xferlogs
usr/local/lib/libguestlib.la
usr/local/lib/libguestlib.so
usr/local/lib/libguestlib.so.0
usr/local/lib/libguestlib.so.0.0.0
usr/local/lib/libhgfs.la
usr/local/lib/libhgfs.so
usr/local/lib/libhgfs.so.0
usr/local/lib/libhgfs.so.0.0.0
usr/local/lib/libvmtools.la
usr/local/lib/libvmtools.so
usr/local/lib/libvmtools.so.0
usr/local/lib/libvmtools.so.0.0.0
usr/local/lib/modules/3.0.21-tinycore/kernel/drivers/misc/vmci.ko
usr/local/lib/modules/3.0.21-tinycore/kernel/drivers/net/vmxnet.ko
usr/local/lib/modules/3.0.21-tinycore/kernel/fs/vmblock/vmblock.ko
usr/local/lib/modules/3.0.21-tinycore/kernel/fs/vmhgfs/vmhgfs.ko
usr/local/lib/modules/3.0.21-tinycore/kernel/net/vsock/vsock.ko
usr/local/lib/open-vm-tools/bin/vmtoolsd
usr/local/lib/open-vm-tools/bin/vmware-checkvm
usr/local/lib/open-vm-tools/bin/vmware-hgfsclient
usr/local/lib/open-vm-tools/bin/vmware-rpctool
usr/local/lib/open-vm-tools/bin/vmware-toolbox-cmd
usr/local/lib/open-vm-tools/bin/vmware-user
usr/local/lib/open-vm-tools/bin/vmware-user-suid-wrapper
usr/local/lib/open-vm-tools/bin/vmware-vmblock-fuse
usr/local/lib/open-vm-tools/bin/vmware-xferlogs
usr/local/lib/open-vm-tools/lib/libguestlib.la
usr/local/lib/open-vm-tools/lib/libguestlib.so
usr/local/lib/open-vm-tools/lib/libguestlib.so.0
usr/local/lib/open-vm-tools/lib/libguestlib.so.0.0.0
usr/local/lib/open-vm-tools/lib/libhgfs.la
usr/local/lib/open-vm-tools/lib/libhgfs.so
usr/local/lib/open-vm-tools/lib/libhgfs.so.0
usr/local/lib/open-vm-tools/lib/libhgfs.so.0.0.0
usr/local/lib/open-vm-tools/lib/libvmtools.la
usr/local/lib/open-vm-tools/lib/libvmtools.so
usr/local/lib/open-vm-tools/lib/libvmtools.so.0
usr/local/lib/open-vm-tools/lib/libvmtools.so.0.0.0
usr/local/lib/open-vm-tools/plugins/common/libhgfsServer.la
usr/local/lib/open-vm-tools/plugins/common/libhgfsServer.so
usr/local/lib/open-vm-tools/plugins/common/libvix.la
usr/local/lib/open-vm-tools/plugins/common/libvix.so
usr/local/lib/open-vm-tools/plugins/vmsvc/libguestInfo.so
usr/local/lib/open-vm-tools/plugins/vmsvc/libpowerOps.so
usr/local/lib/open-vm-tools/plugins/vmsvc/libtimeSync.so
usr/local/lib/open-vm-tools/plugins/vmsvc/libvmbackup.so
usr/local/lib/open-vm-tools/plugins/vmusr/libdesktopEvents.so
usr/local/lib/open-vm-tools/plugins/vmusr/libdndcp.so
usr/local/lib/open-vm-tools/plugins/vmusr/libresolutionSet.so
usr/local/lib/open-vm-tools/sbin/mount.vmhgfs
usr/local/sbin/mount.vmhgfs
usr/local/tce.installed/open-vm-tools-3.0.21-tinycore
