# tinycorelinux/4.x/x86/tcz/postfix.tcz.list
usr/local/bin/mailq
usr/local/bin/newaliases
usr/local/etc/init.d/postfix
usr/local/etc/postfix/access
usr/local/etc/postfix/aliases
usr/local/etc/postfix/bounce.cf.default
usr/local/etc/postfix/canonical
usr/local/etc/postfix/generic
usr/local/etc/postfix/header_checks
usr/local/etc/postfix/main.cf.default
usr/local/etc/postfix/main.cf.sample
usr/local/etc/postfix/makedefs.out
usr/local/etc/postfix/master.cf.sample
usr/local/etc/postfix/relocated
usr/local/etc/postfix/transport
usr/local/etc/postfix/virtual
usr/local/libexec/postfix/anvil
usr/local/libexec/postfix/bounce
usr/local/libexec/postfix/cleanup
usr/local/libexec/postfix/discard
usr/local/libexec/postfix/dnsblog
usr/local/libexec/postfix/error
usr/local/libexec/postfix/flush
usr/local/libexec/postfix/lmtp
usr/local/libexec/postfix/local
usr/local/libexec/postfix/master
usr/local/libexec/postfix/nqmgr
usr/local/libexec/postfix/oqmgr
usr/local/libexec/postfix/pickup
usr/local/libexec/postfix/pipe
usr/local/libexec/postfix/post-install
usr/local/libexec/postfix/postfix-files
usr/local/libexec/postfix/postfix-script
usr/local/libexec/postfix/postfix-wrapper
usr/local/libexec/postfix/postmulti-script
usr/local/libexec/postfix/postscreen
usr/local/libexec/postfix/proxymap
usr/local/libexec/postfix/qmgr
usr/local/libexec/postfix/qmqpd
usr/local/libexec/postfix/scache
usr/local/libexec/postfix/showq
usr/local/libexec/postfix/smtp
usr/local/libexec/postfix/smtpd
usr/local/libexec/postfix/spawn
usr/local/libexec/postfix/tlsmgr
usr/local/libexec/postfix/tlsproxy
usr/local/libexec/postfix/trivial-rewrite
usr/local/libexec/postfix/verify
usr/local/sbin/postalias
usr/local/sbin/postcat
usr/local/sbin/postconf
usr/local/sbin/postdrop
usr/local/sbin/postfix
usr/local/sbin/postkick
usr/local/sbin/postlock
usr/local/sbin/postlog
usr/local/sbin/postmap
usr/local/sbin/postmulti
usr/local/sbin/postqueue
usr/local/sbin/postsuper
usr/local/sbin/sendmail
usr/local/share/doc/License/postfix-copyright.txt
usr/local/share/doc/License/postfix.txt
usr/local/tce.installed/postfix
