# tinycorelinux/4.x/x86/tcz/privoxy.tcz.list
usr/local/etc/init.d/privoxy
usr/local/sbin/privoxy
usr/local/share/doc/License/privoxy-authors.txt
usr/local/share/doc/License/privoxy.txt
usr/local/share/privoxy/conf.tgz
usr/local/tce.installed/privoxy
