# tinycorelinux/4.x/x86/tcz/proftpd.tcz.list
usr/local/bin/ftpasswd
usr/local/bin/ftpcount
usr/local/bin/ftpdctl
usr/local/bin/ftpmail
usr/local/bin/ftpquota
usr/local/bin/ftptop
usr/local/bin/ftpwho
usr/local/bin/prxs
usr/local/etc/init.d/proftpd
usr/local/etc/proftpd.conf
usr/local/lib/pkgconfig/proftpd.pc
usr/local/sbin/ftpasswd
usr/local/sbin/ftpscrub
usr/local/sbin/ftpshut
usr/local/sbin/in.proftpd
usr/local/sbin/proftpd
usr/local/share/locale/bg_BG/LC_MESSAGES/proftpd.mo
usr/local/share/locale/en_US/LC_MESSAGES/proftpd.mo
usr/local/share/locale/fr_FR/LC_MESSAGES/proftpd.mo
usr/local/share/locale/it_IT/LC_MESSAGES/proftpd.mo
usr/local/share/locale/ja_JP/LC_MESSAGES/proftpd.mo
usr/local/share/locale/ko_KR/LC_MESSAGES/proftpd.mo
usr/local/share/locale/ru_RU/LC_MESSAGES/proftpd.mo
usr/local/share/locale/zh_CN/LC_MESSAGES/proftpd.mo
usr/local/share/locale/zh_TW/LC_MESSAGES/proftpd.mo
