# tinycorelinux/4.x/x86/tcz/roam.tcz.list
usr/local/etc/init.d/roam
usr/local/etc/netplug.d/netplug
usr/local/etc/roam.conf.example
usr/local/etc/wpa_supplicant.conf.example
usr/local/sbin/ifdown
usr/local/sbin/ifup
usr/local/tce.installed/roam
