# tinycorelinux/4.x/x86/tcz/rtorrent.tcz.list
usr/local/bin/rtorrent
usr/local/etc/.rtorrent.rc
usr/local/etc/init.d/rtorrent
