# tinycorelinux/4.x/x86/tcz/sixad.tcz.list
usr/local/bin/sixad
usr/local/bin/sixad-jack
usr/local/etc/default/sixad
usr/local/etc/init.d/sixad
usr/local/etc/logrotate.d/sixad
usr/local/sbin/hidraw-dump
usr/local/sbin/sixad-3in1
usr/local/sbin/sixad-bin
usr/local/sbin/sixad-dbus-blocker
usr/local/sbin/sixad-raw
usr/local/sbin/sixad-remote
usr/local/sbin/sixad-sixaxis
usr/local/sbin/sixpair
usr/local/sbin/sixpair-kbd
