# tinycorelinux/4.x/x86/tcz/thelinkbox.tcz.list
usr/local/bin/tlbchat
usr/local/bin/tlbcmd
usr/local/etc/init.d/thelinkbox
usr/local/etc/port.conf.sample
usr/local/etc/tlb.cmds.sample
usr/local/etc/tlb.conf.sample
usr/local/libexec/tlb
