# tinycorelinux/4.x/x86/tcz/tinc.tcz.list
usr/local/etc/init.d/tincd
usr/local/etc/tinc/hosts/alice.sample
usr/local/etc/tinc/howto.txt
usr/local/etc/tinc/tinc-up.sample
usr/local/etc/tinc/tinc.conf.sample
usr/local/sbin/tincd
usr/local/share/doc/License/tinc-authors.txt
usr/local/share/doc/License/tinc-copying.txt
usr/local/share/info/tinc.info
usr/local/share/man/man5/tinc.conf.5.xz
usr/local/share/man/man8/tincd.8.xz
