# tinycorelinux/4.x/x86/tcz/tinyproxy.tcz.list
usr/local/etc/init.d/tinyproxy
usr/local/etc/tinyproxy/tinyproxy.conf
usr/local/sbin/tinyproxy
usr/local/share/man/man5/tinyproxy.conf.5
usr/local/share/man/man8/tinyproxy.8
usr/local/share/tinyproxy/debug.html
usr/local/share/tinyproxy/default.html
usr/local/share/tinyproxy/stats.html
