# tinycorelinux/4.x/x86/tcz/tor.tcz.list
usr/local/bin/tor
usr/local/bin/tor-gencert
usr/local/bin/tor-resolve
usr/local/bin/torify
usr/local/etc/init.d/tor
usr/local/etc/tor/tor-tsocks.conf.sample
usr/local/etc/tor/torrc.sample
usr/local/share/doc/License/tor.txt
usr/local/share/tor/geoip
usr/local/tce.installed/tor
