# tinycorelinux/4.x/x86/tcz/vnstat.tcz.list
usr/local/bin/vnstat
usr/local/etc/init.d/vnstat
usr/local/etc/vnstat.conf
usr/local/sbin/vnstatd
