# tinycorelinux/4.x/x86/tcz/wpa_gui.tcz.list
usr/local/bin/wpa_gui
usr/local/bin/wpa_gui.sh
usr/local/etc/init.d/wpa_gui
usr/local/etc/wpa_gui.conf
usr/local/share/applications/wpa_gui.desktop
usr/local/share/pixmaps/wpa_gui.png
