# tinycorelinux/4.x/x86_64/tcz/alsa.tcz.list
/usr/local/bin/aconnect
/usr/local/bin/alsaloop
/usr/local/bin/alsamixer
/usr/local/bin/alsaucm
/usr/local/bin/amidi
/usr/local/bin/amixer
/usr/local/bin/aplay
/usr/local/bin/aplaymidi
/usr/local/bin/arecord
/usr/local/bin/arecordmidi
/usr/local/bin/aseqdump
/usr/local/bin/aseqnet
/usr/local/bin/aserver
/usr/local/bin/iecset
/usr/local/bin/speaker-test
/usr/local/etc/init.d/alsasound
/usr/local/lib/alsa-lib/smixer/smixer-ac97.so
/usr/local/lib/alsa-lib/smixer/smixer-hda.so
/usr/local/lib/alsa-lib/smixer/smixer-sbase.so
/usr/local/sbin/alsactl
/usr/local/share/sounds/alsa/Front_Center.wav
/usr/local/share/sounds/alsa/Front_Left.wav
/usr/local/share/sounds/alsa/Front_Right.wav
/usr/local/share/sounds/alsa/Noise.wav
/usr/local/share/sounds/alsa/Rear_Center.wav
/usr/local/share/sounds/alsa/Rear_Left.wav
/usr/local/share/sounds/alsa/Rear_Right.wav
/usr/local/share/sounds/alsa/Side_Left.wav
/usr/local/share/sounds/alsa/Side_Right.wav
