# tinycorelinux/15.x/armhf/tcz/mesa-dev.tcz.list
/usr/local/include/GL/gl.h
/usr/local/include/GL/glcorearb.h
/usr/local/include/GL/glext.h
/usr/local/include/GL/glx.h
/usr/local/include/GL/glxext.h
/usr/local/include/GL/internal/dri_interface.h
/usr/local/include/KHR/khrplatform.h
/usr/local/lib/pkgconfig/dri.pc
/usr/local/lib/pkgconfig/gl.pc
/usr/local/share/drirc.d/00-mesa-defaults.conf
/usr/local/share/vulkan/icd.d/broadcom_icd.armv7l.json
/usr/local/share/vulkan/icd.d/lvp_icd.armv7l.json
