# tinycorelinux/8.x/x86/tcz/apache2.4.tcz.list
/usr/local/apache2/cgi-bin/printenv
/usr/local/apache2/cgi-bin/printenv.vbs
/usr/local/apache2/cgi-bin/printenv.wsf
/usr/local/apache2/cgi-bin/test-cgi
/usr/local/apache2/error/HTTP_BAD_GATEWAY.html.var
/usr/local/apache2/error/HTTP_BAD_REQUEST.html.var
/usr/local/apache2/error/HTTP_FORBIDDEN.html.var
/usr/local/apache2/error/HTTP_GONE.html.var
/usr/local/apache2/error/HTTP_INTERNAL_SERVER_ERROR.html.var
/usr/local/apache2/error/HTTP_LENGTH_REQUIRED.html.var
/usr/local/apache2/error/HTTP_METHOD_NOT_ALLOWED.html.var
/usr/local/apache2/error/HTTP_NOT_FOUND.html.var
/usr/local/apache2/error/HTTP_NOT_IMPLEMENTED.html.var
/usr/local/apache2/error/HTTP_PRECONDITION_FAILED.html.var
/usr/local/apache2/error/HTTP_REQUEST_ENTITY_TOO_LARGE.html.var
/usr/local/apache2/error/HTTP_REQUEST_TIME_OUT.html.var
/usr/local/apache2/error/HTTP_REQUEST_URI_TOO_LARGE.html.var
/usr/local/apache2/error/HTTP_SERVICE_UNAVAILABLE.html.var
/usr/local/apache2/error/HTTP_UNAUTHORIZED.html.var
/usr/local/apache2/error/HTTP_UNSUPPORTED_MEDIA_TYPE.html.var
/usr/local/apache2/error/HTTP_VARIANT_ALSO_VARIES.html.var
/usr/local/apache2/error/README
/usr/local/apache2/error/contact.html.var
/usr/local/apache2/error/include/bottom.html
/usr/local/apache2/error/include/spacer.html
/usr/local/apache2/error/include/top.html
/usr/local/apache2/htdocs/index.html
/usr/local/apache2/icons/README
/usr/local/apache2/icons/README.html
/usr/local/apache2/icons/a.gif
/usr/local/apache2/icons/a.png
/usr/local/apache2/icons/alert.black.gif
/usr/local/apache2/icons/alert.black.png
/usr/local/apache2/icons/alert.red.gif
/usr/local/apache2/icons/alert.red.png
/usr/local/apache2/icons/apache_pb.gif
/usr/local/apache2/icons/apache_pb.png
/usr/local/apache2/icons/apache_pb.svg
/usr/local/apache2/icons/apache_pb2.gif
/usr/local/apache2/icons/apache_pb2.png
/usr/local/apache2/icons/back.gif
/usr/local/apache2/icons/back.png
/usr/local/apache2/icons/ball.gray.gif
/usr/local/apache2/icons/ball.gray.png
/usr/local/apache2/icons/ball.red.gif
/usr/local/apache2/icons/ball.red.png
/usr/local/apache2/icons/binary.gif
/usr/local/apache2/icons/binary.png
/usr/local/apache2/icons/binhex.gif
/usr/local/apache2/icons/binhex.png
/usr/local/apache2/icons/blank.gif
/usr/local/apache2/icons/blank.png
/usr/local/apache2/icons/bomb.gif
/usr/local/apache2/icons/bomb.png
/usr/local/apache2/icons/box1.gif
/usr/local/apache2/icons/box1.png
/usr/local/apache2/icons/box2.gif
/usr/local/apache2/icons/box2.png
/usr/local/apache2/icons/broken.gif
/usr/local/apache2/icons/broken.png
/usr/local/apache2/icons/burst.gif
/usr/local/apache2/icons/burst.png
/usr/local/apache2/icons/c.gif
/usr/local/apache2/icons/c.png
/usr/local/apache2/icons/comp.blue.gif
/usr/local/apache2/icons/comp.blue.png
/usr/local/apache2/icons/comp.gray.gif
/usr/local/apache2/icons/comp.gray.png
/usr/local/apache2/icons/compressed.gif
/usr/local/apache2/icons/compressed.png
/usr/local/apache2/icons/continued.gif
/usr/local/apache2/icons/continued.png
/usr/local/apache2/icons/dir.gif
/usr/local/apache2/icons/dir.png
/usr/local/apache2/icons/diskimg.gif
/usr/local/apache2/icons/diskimg.png
/usr/local/apache2/icons/down.gif
/usr/local/apache2/icons/down.png
/usr/local/apache2/icons/dvi.gif
/usr/local/apache2/icons/dvi.png
/usr/local/apache2/icons/f.gif
/usr/local/apache2/icons/f.png
/usr/local/apache2/icons/folder.gif
/usr/local/apache2/icons/folder.open.gif
/usr/local/apache2/icons/folder.open.png
/usr/local/apache2/icons/folder.png
/usr/local/apache2/icons/folder.sec.gif
/usr/local/apache2/icons/folder.sec.png
/usr/local/apache2/icons/forward.gif
/usr/local/apache2/icons/forward.png
/usr/local/apache2/icons/generic.gif
/usr/local/apache2/icons/generic.png
/usr/local/apache2/icons/generic.red.gif
/usr/local/apache2/icons/generic.red.png
/usr/local/apache2/icons/generic.sec.gif
/usr/local/apache2/icons/generic.sec.png
/usr/local/apache2/icons/hand.right.gif
/usr/local/apache2/icons/hand.right.png
/usr/local/apache2/icons/hand.up.gif
/usr/local/apache2/icons/hand.up.png
/usr/local/apache2/icons/icon.sheet.gif
/usr/local/apache2/icons/icon.sheet.png
/usr/local/apache2/icons/image1.gif
/usr/local/apache2/icons/image1.png
/usr/local/apache2/icons/image2.gif
/usr/local/apache2/icons/image2.png
/usr/local/apache2/icons/image3.gif
/usr/local/apache2/icons/image3.png
/usr/local/apache2/icons/index.gif
/usr/local/apache2/icons/index.png
/usr/local/apache2/icons/layout.gif
/usr/local/apache2/icons/layout.png
/usr/local/apache2/icons/left.gif
/usr/local/apache2/icons/left.png
/usr/local/apache2/icons/link.gif
/usr/local/apache2/icons/link.png
/usr/local/apache2/icons/movie.gif
/usr/local/apache2/icons/movie.png
/usr/local/apache2/icons/odf6odb.png
/usr/local/apache2/icons/odf6odc.png
/usr/local/apache2/icons/odf6odf.png
/usr/local/apache2/icons/odf6odg.png
/usr/local/apache2/icons/odf6odi.png
/usr/local/apache2/icons/odf6odm.png
/usr/local/apache2/icons/odf6odp.png
/usr/local/apache2/icons/odf6ods.png
/usr/local/apache2/icons/odf6odt.png
/usr/local/apache2/icons/odf6otc.png
/usr/local/apache2/icons/odf6otf.png
/usr/local/apache2/icons/odf6otg.png
/usr/local/apache2/icons/odf6oth.png
/usr/local/apache2/icons/odf6oti.png
/usr/local/apache2/icons/odf6otp.png
/usr/local/apache2/icons/odf6ots.png
/usr/local/apache2/icons/odf6ott.png
/usr/local/apache2/icons/p.gif
/usr/local/apache2/icons/p.png
/usr/local/apache2/icons/patch.gif
/usr/local/apache2/icons/patch.png
/usr/local/apache2/icons/pdf.gif
/usr/local/apache2/icons/pdf.png
/usr/local/apache2/icons/pie0.gif
/usr/local/apache2/icons/pie0.png
/usr/local/apache2/icons/pie1.gif
/usr/local/apache2/icons/pie1.png
/usr/local/apache2/icons/pie2.gif
/usr/local/apache2/icons/pie2.png
/usr/local/apache2/icons/pie3.gif
/usr/local/apache2/icons/pie3.png
/usr/local/apache2/icons/pie4.gif
/usr/local/apache2/icons/pie4.png
/usr/local/apache2/icons/pie5.gif
/usr/local/apache2/icons/pie5.png
/usr/local/apache2/icons/pie6.gif
/usr/local/apache2/icons/pie6.png
/usr/local/apache2/icons/pie7.gif
/usr/local/apache2/icons/pie7.png
/usr/local/apache2/icons/pie8.gif
/usr/local/apache2/icons/pie8.png
/usr/local/apache2/icons/portal.gif
/usr/local/apache2/icons/portal.png
/usr/local/apache2/icons/ps.gif
/usr/local/apache2/icons/ps.png
/usr/local/apache2/icons/quill.gif
/usr/local/apache2/icons/quill.png
/usr/local/apache2/icons/right.gif
/usr/local/apache2/icons/right.png
/usr/local/apache2/icons/screw1.gif
/usr/local/apache2/icons/screw1.png
/usr/local/apache2/icons/screw2.gif
/usr/local/apache2/icons/screw2.png
/usr/local/apache2/icons/script.gif
/usr/local/apache2/icons/script.png
/usr/local/apache2/icons/small/back.gif
/usr/local/apache2/icons/small/back.png
/usr/local/apache2/icons/small/binary.gif
/usr/local/apache2/icons/small/binary.png
/usr/local/apache2/icons/small/binhex.gif
/usr/local/apache2/icons/small/binhex.png
/usr/local/apache2/icons/small/blank.gif
/usr/local/apache2/icons/small/blank.png
/usr/local/apache2/icons/small/broken.gif
/usr/local/apache2/icons/small/broken.png
/usr/local/apache2/icons/small/burst.gif
/usr/local/apache2/icons/small/burst.png
/usr/local/apache2/icons/small/comp1.gif
/usr/local/apache2/icons/small/comp1.png
/usr/local/apache2/icons/small/comp2.gif
/usr/local/apache2/icons/small/comp2.png
/usr/local/apache2/icons/small/compressed.gif
/usr/local/apache2/icons/small/compressed.png
/usr/local/apache2/icons/small/continued.gif
/usr/local/apache2/icons/small/continued.png
/usr/local/apache2/icons/small/doc.gif
/usr/local/apache2/icons/small/doc.png
/usr/local/apache2/icons/small/folder.gif
/usr/local/apache2/icons/small/folder.png
/usr/local/apache2/icons/small/folder2.gif
/usr/local/apache2/icons/small/folder2.png
/usr/local/apache2/icons/small/forward.gif
/usr/local/apache2/icons/small/forward.png
/usr/local/apache2/icons/small/generic.gif
/usr/local/apache2/icons/small/generic.png
/usr/local/apache2/icons/small/generic2.gif
/usr/local/apache2/icons/small/generic2.png
/usr/local/apache2/icons/small/generic3.gif
/usr/local/apache2/icons/small/generic3.png
/usr/local/apache2/icons/small/image.gif
/usr/local/apache2/icons/small/image.png
/usr/local/apache2/icons/small/image2.gif
/usr/local/apache2/icons/small/image2.png
/usr/local/apache2/icons/small/index.gif
/usr/local/apache2/icons/small/index.png
/usr/local/apache2/icons/small/key.gif
/usr/local/apache2/icons/small/key.png
/usr/local/apache2/icons/small/movie.gif
/usr/local/apache2/icons/small/movie.png
/usr/local/apache2/icons/small/patch.gif
/usr/local/apache2/icons/small/patch.png
/usr/local/apache2/icons/small/ps.gif
/usr/local/apache2/icons/small/ps.png
/usr/local/apache2/icons/small/rainbow.gif
/usr/local/apache2/icons/small/rainbow.png
/usr/local/apache2/icons/small/sound.gif
/usr/local/apache2/icons/small/sound.png
/usr/local/apache2/icons/small/sound2.gif
/usr/local/apache2/icons/small/sound2.png
/usr/local/apache2/icons/small/tar.gif
/usr/local/apache2/icons/small/tar.png
/usr/local/apache2/icons/small/text.gif
/usr/local/apache2/icons/small/text.png
/usr/local/apache2/icons/small/transfer.gif
/usr/local/apache2/icons/small/transfer.png
/usr/local/apache2/icons/small/unknown.gif
/usr/local/apache2/icons/small/unknown.png
/usr/local/apache2/icons/small/uu.gif
/usr/local/apache2/icons/small/uu.png
/usr/local/apache2/icons/sound1.gif
/usr/local/apache2/icons/sound1.png
/usr/local/apache2/icons/sound2.gif
/usr/local/apache2/icons/sound2.png
/usr/local/apache2/icons/sphere1.gif
/usr/local/apache2/icons/sphere1.png
/usr/local/apache2/icons/sphere2.gif
/usr/local/apache2/icons/sphere2.png
/usr/local/apache2/icons/svg.png
/usr/local/apache2/icons/tar.gif
/usr/local/apache2/icons/tar.png
/usr/local/apache2/icons/tex.gif
/usr/local/apache2/icons/tex.png
/usr/local/apache2/icons/text.gif
/usr/local/apache2/icons/text.png
/usr/local/apache2/icons/transfer.gif
/usr/local/apache2/icons/transfer.png
/usr/local/apache2/icons/unknown.gif
/usr/local/apache2/icons/unknown.png
/usr/local/apache2/icons/up.gif
/usr/local/apache2/icons/up.png
/usr/local/apache2/icons/uu.gif
/usr/local/apache2/icons/uu.png
/usr/local/apache2/icons/uuencoded.gif
/usr/local/apache2/icons/uuencoded.png
/usr/local/apache2/icons/world1.gif
/usr/local/apache2/icons/world1.png
/usr/local/apache2/icons/world2.gif
/usr/local/apache2/icons/world2.png
/usr/local/apache2/icons/xml.png
/usr/local/apache2/modules/httpd.exp
/usr/local/apache2/modules/mod_access_compat.so
/usr/local/apache2/modules/mod_actions.so
/usr/local/apache2/modules/mod_alias.so
/usr/local/apache2/modules/mod_allowmethods.so
/usr/local/apache2/modules/mod_asis.so
/usr/local/apache2/modules/mod_auth_basic.so
/usr/local/apache2/modules/mod_auth_digest.so
/usr/local/apache2/modules/mod_auth_form.so
/usr/local/apache2/modules/mod_authn_anon.so
/usr/local/apache2/modules/mod_authn_core.so
/usr/local/apache2/modules/mod_authn_dbd.so
/usr/local/apache2/modules/mod_authn_dbm.so
/usr/local/apache2/modules/mod_authn_file.so
/usr/local/apache2/modules/mod_authn_socache.so
/usr/local/apache2/modules/mod_authnz_fcgi.so
/usr/local/apache2/modules/mod_authnz_ldap.so
/usr/local/apache2/modules/mod_authz_core.so
/usr/local/apache2/modules/mod_authz_dbd.so
/usr/local/apache2/modules/mod_authz_dbm.so
/usr/local/apache2/modules/mod_authz_groupfile.so
/usr/local/apache2/modules/mod_authz_host.so
/usr/local/apache2/modules/mod_authz_owner.so
/usr/local/apache2/modules/mod_authz_user.so
/usr/local/apache2/modules/mod_autoindex.so
/usr/local/apache2/modules/mod_bucketeer.so
/usr/local/apache2/modules/mod_buffer.so
/usr/local/apache2/modules/mod_cache.so
/usr/local/apache2/modules/mod_cache_disk.so
/usr/local/apache2/modules/mod_cache_socache.so
/usr/local/apache2/modules/mod_case_filter.so
/usr/local/apache2/modules/mod_case_filter_in.so
/usr/local/apache2/modules/mod_cern_meta.so
/usr/local/apache2/modules/mod_cgi.so
/usr/local/apache2/modules/mod_cgid.so
/usr/local/apache2/modules/mod_charset_lite.so
/usr/local/apache2/modules/mod_data.so
/usr/local/apache2/modules/mod_dav.so
/usr/local/apache2/modules/mod_dav_fs.so
/usr/local/apache2/modules/mod_dav_lock.so
/usr/local/apache2/modules/mod_dbd.so
/usr/local/apache2/modules/mod_deflate.so
/usr/local/apache2/modules/mod_dialup.so
/usr/local/apache2/modules/mod_dir.so
/usr/local/apache2/modules/mod_dumpio.so
/usr/local/apache2/modules/mod_echo.so
/usr/local/apache2/modules/mod_env.so
/usr/local/apache2/modules/mod_example_hooks.so
/usr/local/apache2/modules/mod_example_ipc.so
/usr/local/apache2/modules/mod_expires.so
/usr/local/apache2/modules/mod_ext_filter.so
/usr/local/apache2/modules/mod_file_cache.so
/usr/local/apache2/modules/mod_filter.so
/usr/local/apache2/modules/mod_headers.so
/usr/local/apache2/modules/mod_heartbeat.so
/usr/local/apache2/modules/mod_heartmonitor.so
/usr/local/apache2/modules/mod_http2.so
/usr/local/apache2/modules/mod_ident.so
/usr/local/apache2/modules/mod_imagemap.so
/usr/local/apache2/modules/mod_include.so
/usr/local/apache2/modules/mod_info.so
/usr/local/apache2/modules/mod_isapi.so
/usr/local/apache2/modules/mod_lbmethod_bybusyness.so
/usr/local/apache2/modules/mod_lbmethod_byrequests.so
/usr/local/apache2/modules/mod_lbmethod_bytraffic.so
/usr/local/apache2/modules/mod_lbmethod_heartbeat.so
/usr/local/apache2/modules/mod_ldap.so
/usr/local/apache2/modules/mod_log_config.so
/usr/local/apache2/modules/mod_log_debug.so
/usr/local/apache2/modules/mod_log_forensic.so
/usr/local/apache2/modules/mod_logio.so
/usr/local/apache2/modules/mod_lua.so
/usr/local/apache2/modules/mod_macro.so
/usr/local/apache2/modules/mod_mime.so
/usr/local/apache2/modules/mod_mime_magic.so
/usr/local/apache2/modules/mod_mpm_event.so
/usr/local/apache2/modules/mod_mpm_prefork.so
/usr/local/apache2/modules/mod_mpm_worker.so
/usr/local/apache2/modules/mod_negotiation.so
/usr/local/apache2/modules/mod_optional_fn_export.so
/usr/local/apache2/modules/mod_optional_fn_import.so
/usr/local/apache2/modules/mod_optional_hook_export.so
/usr/local/apache2/modules/mod_optional_hook_import.so
/usr/local/apache2/modules/mod_proxy.so
/usr/local/apache2/modules/mod_proxy_ajp.so
/usr/local/apache2/modules/mod_proxy_balancer.so
/usr/local/apache2/modules/mod_proxy_connect.so
/usr/local/apache2/modules/mod_proxy_express.so
/usr/local/apache2/modules/mod_proxy_fcgi.so
/usr/local/apache2/modules/mod_proxy_fdpass.so
/usr/local/apache2/modules/mod_proxy_ftp.so
/usr/local/apache2/modules/mod_proxy_hcheck.so
/usr/local/apache2/modules/mod_proxy_html.so
/usr/local/apache2/modules/mod_proxy_http.so
/usr/local/apache2/modules/mod_proxy_http2.so
/usr/local/apache2/modules/mod_proxy_scgi.so
/usr/local/apache2/modules/mod_proxy_wstunnel.so
/usr/local/apache2/modules/mod_ratelimit.so
/usr/local/apache2/modules/mod_reflector.so
/usr/local/apache2/modules/mod_remoteip.so
/usr/local/apache2/modules/mod_reqtimeout.so
/usr/local/apache2/modules/mod_request.so
/usr/local/apache2/modules/mod_rewrite.so
/usr/local/apache2/modules/mod_sed.so
/usr/local/apache2/modules/mod_session.so
/usr/local/apache2/modules/mod_session_cookie.so
/usr/local/apache2/modules/mod_session_crypto.so
/usr/local/apache2/modules/mod_session_dbd.so
/usr/local/apache2/modules/mod_setenvif.so
/usr/local/apache2/modules/mod_slotmem_plain.so
/usr/local/apache2/modules/mod_slotmem_shm.so
/usr/local/apache2/modules/mod_socache_dbm.so
/usr/local/apache2/modules/mod_socache_memcache.so
/usr/local/apache2/modules/mod_socache_shmcb.so
/usr/local/apache2/modules/mod_speling.so
/usr/local/apache2/modules/mod_ssl.so
/usr/local/apache2/modules/mod_status.so
/usr/local/apache2/modules/mod_substitute.so
/usr/local/apache2/modules/mod_suexec.so
/usr/local/apache2/modules/mod_unique_id.so
/usr/local/apache2/modules/mod_unixd.so
/usr/local/apache2/modules/mod_userdir.so
/usr/local/apache2/modules/mod_usertrack.so
/usr/local/apache2/modules/mod_version.so
/usr/local/apache2/modules/mod_vhost_alias.so
/usr/local/apache2/modules/mod_watchdog.so
/usr/local/apache2/modules/mod_xml2enc.so
/usr/local/bin/ab
/usr/local/bin/dbmmanage
/usr/local/bin/htdbm
/usr/local/bin/htdigest
/usr/local/bin/htpasswd
/usr/local/bin/httxt2dbm
/usr/local/bin/logresolve
/usr/local/etc/httpd/magic
/usr/local/etc/httpd/mime.types
/usr/local/etc/httpd/original/conf.d/httpd-autoindex.conf
/usr/local/etc/httpd/original/conf.d/httpd-dav.conf
/usr/local/etc/httpd/original/conf.d/httpd-default.conf
/usr/local/etc/httpd/original/conf.d/httpd-info.conf
/usr/local/etc/httpd/original/conf.d/httpd-languages.conf
/usr/local/etc/httpd/original/conf.d/httpd-manual.conf
/usr/local/etc/httpd/original/conf.d/httpd-mpm.conf
/usr/local/etc/httpd/original/conf.d/httpd-multilang-errordoc.conf
/usr/local/etc/httpd/original/conf.d/httpd-ssl.conf
/usr/local/etc/httpd/original/conf.d/httpd-userdir.conf
/usr/local/etc/httpd/original/conf.d/httpd-vhosts.conf
/usr/local/etc/httpd/original/conf.d/proxy-html.conf
/usr/local/etc/httpd/original/httpd.conf-sample
/usr/local/etc/init.d/httpd
/usr/local/sbin/apachectl
/usr/local/sbin/checkgid
/usr/local/sbin/envvars
/usr/local/sbin/envvars-std
/usr/local/sbin/fcgistarter
/usr/local/sbin/htcacheclean
/usr/local/sbin/httpd
/usr/local/sbin/rotatelogs
/usr/local/sbin/suexec
