# tinycorelinux/8.x/x86/tcz/avahi.tcz.list
/usr/local/bin/avahi-browse
/usr/local/bin/avahi-browse-domains
/usr/local/bin/avahi-publish
/usr/local/bin/avahi-publish-address
/usr/local/bin/avahi-publish-service
/usr/local/bin/avahi-resolve
/usr/local/bin/avahi-resolve-address
/usr/local/bin/avahi-resolve-host-name
/usr/local/bin/avahi-set-host-name
/usr/local/etc/avahi/avahi-autoipd.action
/usr/local/etc/avahi/avahi-dnsconfd.action
/usr/local/etc/avahi/services/sftp-ssh.service
/usr/local/etc/avahi/services/ssh.service
/usr/local/etc/dbus-1/system.d/avahi-dbus.conf
/usr/local/etc/init.d/avahi
/usr/local/sbin/avahi-autoipd
/usr/local/sbin/avahi-daemon
/usr/local/sbin/avahi-dnsconfd
/usr/local/share/avahi/avahi-service.dtd
/usr/local/share/avahi/files
/usr/local/share/avahi/files/avahi-daemon.conf
/usr/local/share/avahi/files/hosts
/usr/local/share/avahi/service-types
/usr/local/share/dbus-1/interfaces/org.freedesktop.Avahi.AddressResolver.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.Avahi.DomainBrowser.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.Avahi.EntryGroup.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.Avahi.HostNameResolver.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.Avahi.RecordBrowser.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.Avahi.Server.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.Avahi.ServiceBrowser.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.Avahi.ServiceResolver.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.Avahi.ServiceTypeBrowser.xml
/usr/local/tce.installed/avahi
