# tinycorelinux/6.x/x86_64/tcz/rfkill.tcz.list
/local/etc/conf.d/rfkill
/local/etc/init.d/rfkill
/local/sbin/rfkill
/local/share/rfkill/COPYING
