# tinycorelinux/15.x/x86_64/tcz/ifplugd.tcz.list
/usr/local/etc/ifplugd/ifplugd.action.sample
/usr/local/etc/ifplugd/ifplugd.conf.sample
/usr/local/etc/init.d/ifplugd
/usr/local/share/doc/ifplugd/LICENSE
/usr/local/sbin/ifplugstatus
/usr/local/sbin/ifplugd
