# tinycorelinux/15.x/x86_64/tcz/portmap.tcz.list
/usr/local/etc/init.d/portmap
/usr/local/sbin/pmap_dump
/usr/local/sbin/pmap_set
/usr/local/sbin/portmap
