# tinycorelinux/13.x/x86_64/tcz/gpm.tcz.list
/usr/local/bin/disable-paste
/usr/local/bin/display-buttons
/usr/local/bin/display-coords
/usr/local/bin/get-versions
/usr/local/bin/gpm-root
/usr/local/bin/hltest
/usr/local/bin/mev
/usr/local/bin/mouse-test
/usr/local/etc/init.d/gpm
/usr/local/include/gpm.h
/usr/local/lib/libgpm.so
/usr/local/lib/libgpm.so.2
/usr/local/lib/libgpm.so.2.1.0
/usr/local/sbin/gpm
/usr/local/share/gpm/files/gpm-root.conf
/usr/local/share/gpm/files/gpm.conf
/usr/local/tce.installed/gpm
