# tinycorelinux/13.x/x86_64/tcz/redis.tcz.list
/usr/local/bin/redis-benchmark
/usr/local/bin/redis-cli
/usr/local/bin/redis-server
/usr/local/etc/init.d/redis
/usr/local/etc/redis/redis.conf.default
/usr/local/share/doc/redis-6.2.6/COPYING
/usr/local/share/doc/redis-6.2.6/copyright
/usr/local/tce.installed/redis
