# tinycorelinux/4.x/armv6/tcz/dbus.tcz.list
usr/local/tce.installed/dbus.ins
usr/local/etc/init.d/dbus
usr/local/etc/dbus-1/system.conf
usr/local/etc/dbus-1/session.conf
usr/local/libexec/dbus-daemon-launch-helper
usr/local/bin/dbus-daemon
usr/local/bin/dbus-launch
usr/local/bin/dbus-monitor
usr/local/bin/dbus-send
usr/local/bin/dbus-cleanup-sockets
usr/local/bin/dbus-uuidgen
usr/local/lib/libdbus-1.so.3.7.2
usr/local/lib/dbus-1.0/include/dbus/dbus-arch-deps.h
usr/local/lib/libdbus-1.so
usr/local/lib/libdbus-1.so.3
