# tinycorelinux/4.x/armv6/tcz/dropbear.tcz.list
usr/local/bin/dbclient
usr/local/bin/dropbear
usr/local/bin/dropbearconvert
usr/local/bin/dropbearkey
usr/local/bin/scp
usr/local/bin/ssh
usr/local/etc/dropbear/banner
usr/local/etc/init.d/dropbear
usr/local/sbin/dropbearmulti
