# tinycorelinux/11.x/x86_64/tcz/networkmanager.tcz.list
/usr/local/bin/nm-online
/usr/local/bin/nmcli
/usr/local/etc/NetworkManager/nm-system-settings.conf
/usr/local/etc/dbus-1/system.d/nm-dispatcher.conf
/usr/local/etc/dbus-1/system.d/org.freedesktop.NetworkManager.conf
/usr/local/etc/default/NetworkManager
/usr/local/etc/init.d/networkmanager
/usr/local/etc/polkit-1/rules.d/50-org.freedesktop.NetworkManager.rules
/usr/local/lib/NetworkManager/1.22.8/libnm-device-plugin-adsl.so
/usr/local/lib/NetworkManager/1.22.8/libnm-device-plugin-bluetooth.so
/usr/local/lib/NetworkManager/1.22.8/libnm-device-plugin-ovs.so
/usr/local/lib/NetworkManager/1.22.8/libnm-device-plugin-wifi.so
/usr/local/lib/NetworkManager/1.22.8/libnm-device-plugin-wwan.so
/usr/local/lib/NetworkManager/1.22.8/libnm-ppp-plugin.so
/usr/local/lib/NetworkManager/1.22.8/libnm-wwan.so
/usr/local/lib/NetworkManager/nm-cloud-setup
/usr/local/lib/NetworkManager/nm-dhcp-helper
/usr/local/lib/NetworkManager/nm-dispatcher
/usr/local/lib/NetworkManager/nm-iface-helper
/usr/local/lib/NetworkManager/nm-initrd-generator
/usr/local/lib/libnm.so
/usr/local/lib/libnm.so.0
/usr/local/lib/libnm.so.0.1.0
/usr/local/lib/pppd/2.4.5/nm-pppd-plugin.so
/usr/local/sbin/NetworkManager
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.AccessPoint.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.AgentManager.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.Checkpoint.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.Connection.Active.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.DHCP4Config.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.DHCP6Config.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.Device.Adsl.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.Device.Bluetooth.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.Device.Bond.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.Device.Bridge.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.Device.Dummy.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.Device.Generic.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.Device.IPTunnel.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.Device.Infiniband.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.Device.Lowpan.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.Device.Macsec.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.Device.Macvlan.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.Device.Modem.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.Device.OlpcMesh.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.Device.OvsBridge.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.Device.OvsInterface.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.Device.OvsPort.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.Device.Ppp.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.Device.Statistics.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.Device.Team.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.Device.Tun.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.Device.Veth.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.Device.Vlan.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.Device.Vxlan.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.Device.WiMax.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.Device.WifiP2P.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.Device.WireGuard.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.Device.Wired.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.Device.Wireless.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.Device.Wpan.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.Device.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.DnsManager.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.IP4Config.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.IP6Config.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.PPP.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.SecretAgent.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.Settings.Connection.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.Settings.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.VPN.Connection.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.VPN.Plugin.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.WiMax.Nsp.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.WifiP2PPeer.xml
/usr/local/share/dbus-1/interfaces/org.freedesktop.NetworkManager.xml
/usr/local/share/dbus-1/system-services/org.freedesktop.nm_dispatcher.service
/usr/local/share/networkmanager/files/84-nm-drivers.rules
/usr/local/share/networkmanager/files/85-nm-unmanaged.rules
/usr/local/share/networkmanager/files/90-nm-thunderbolt.rules
/usr/local/share/networkmanager/files/server.conf
/usr/local/share/polkit-1/actions/org.freedesktop.NetworkManager.policy
/usr/local/tce.installed/networkmanager
