# tinycorelinux/11.x/x86_64/tcz/nfs-utils.tcz.list
/usr/local/etc/init.d/nfs-client
/usr/local/etc/init.d/nfs-server
/usr/local/sbin/exportfs
/usr/local/sbin/mount.nfs
/usr/local/sbin/mountstats
/usr/local/sbin/nfsiostat
/usr/local/sbin/nfsstat
/usr/local/sbin/osd_login
/usr/local/sbin/rpc.mountd
/usr/local/sbin/rpc.nfsd
/usr/local/sbin/rpc.statd
/usr/local/sbin/rpcdebug
/usr/local/sbin/showmount
/usr/local/sbin/sm-notify
/usr/local/sbin/start-statd
/usr/local/sbin/umount.nfs
/usr/local/share/applications/nfs-utils.desktop
/usr/local/share/applications/nfs-utils~1.desktop
/usr/local/share/applications/nfs-utils~2.desktop
/usr/local/share/applications/nfs-utils~3.desktop
/usr/local/share/applications/nfs-utils~4.desktop
/usr/local/share/applications/nfs-utils~5.desktop
/usr/local/share/nfs-utils/files/exports
/usr/local/tce.installed/nfs-utils
