# tinycorelinux/15.x/aarch64/tcz/linux-pam.tcz.list
/usr/local/etc/environment
/usr/local/etc/pam.d/other
/usr/local/etc/pam.d/system-auth
/usr/local/etc/security/access.conf
/usr/local/etc/security/faillock.conf
/usr/local/etc/security/group.conf
/usr/local/etc/security/limits.conf
/usr/local/etc/security/namespace.conf
/usr/local/etc/security/namespace.init
/usr/local/etc/security/pam_env.conf
/usr/local/etc/security/time.conf
/usr/local/lib/libpam.so
/usr/local/lib/libpam.so.0
/usr/local/lib/libpam.so.0.85.1
/usr/local/lib/libpam_misc.so
/usr/local/lib/libpam_misc.so.0
/usr/local/lib/libpam_misc.so.0.82.1
/usr/local/lib/libpamc.so
/usr/local/lib/libpamc.so.0
/usr/local/lib/libpamc.so.0.82.1
/usr/local/lib/security/pam_access.so
/usr/local/lib/security/pam_debug.so
/usr/local/lib/security/pam_deny.so
/usr/local/lib/security/pam_echo.so
/usr/local/lib/security/pam_env.so
/usr/local/lib/security/pam_exec.so
/usr/local/lib/security/pam_faildelay.so
/usr/local/lib/security/pam_faillock.so
/usr/local/lib/security/pam_filter/upperLOWER
/usr/local/lib/security/pam_filter.so
/usr/local/lib/security/pam_ftp.so
/usr/local/lib/security/pam_group.so
/usr/local/lib/security/pam_issue.so
/usr/local/lib/security/pam_keyinit.so
/usr/local/lib/security/pam_lastlog.so
/usr/local/lib/security/pam_limits.so
/usr/local/lib/security/pam_listfile.so
/usr/local/lib/security/pam_localuser.so
/usr/local/lib/security/pam_loginuid.so
/usr/local/lib/security/pam_mail.so
/usr/local/lib/security/pam_mkhomedir.so
/usr/local/lib/security/pam_motd.so
/usr/local/lib/security/pam_namespace.so
/usr/local/lib/security/pam_nologin.so
/usr/local/lib/security/pam_permit.so
/usr/local/lib/security/pam_pwhistory.so
/usr/local/lib/security/pam_rhosts.so
/usr/local/lib/security/pam_rootok.so
/usr/local/lib/security/pam_securetty.so
/usr/local/lib/security/pam_setquota.so
/usr/local/lib/security/pam_shells.so
/usr/local/lib/security/pam_stress.so
/usr/local/lib/security/pam_succeed_if.so
/usr/local/lib/security/pam_time.so
/usr/local/lib/security/pam_timestamp.so
/usr/local/lib/security/pam_umask.so
/usr/local/lib/security/pam_unix.so
/usr/local/lib/security/pam_usertype.so
/usr/local/lib/security/pam_warn.so
/usr/local/lib/security/pam_wheel.so
/usr/local/lib/security/pam_xauth.so
/usr/local/lib/systemd/system/pam_namespace.service
/usr/local/sbin/faillock
/usr/local/sbin/mkhomedir_helper
/usr/local/sbin/pam_namespace_helper
/usr/local/sbin/pam_timestamp_check
/usr/local/sbin/unix_chkpwd
/usr/local/sbin/unix_update
/usr/local/tce.installed/linux-pam
