# tinycorelinux/16.x/x86/tcz/cpufreqd.tcz.list
/usr/local/bin/cpufreqd-get
/usr/local/bin/cpufreqd-set
/usr/local/etc/init.d/cpufreqd
/usr/local/lib/cpufreqd_acpi.so
/usr/local/lib/cpufreqd_apm.so
/usr/local/lib/cpufreqd_cpu.so
/usr/local/lib/cpufreqd_exec.so
/usr/local/lib/cpufreqd_governor_parameters.so
/usr/local/lib/cpufreqd_nforce2.so
/usr/local/lib/cpufreqd_pmu.so
/usr/local/lib/cpufreqd_programs.so
/usr/local/lib/cpufreqd_tau.so
/usr/local/sbin/cpufreqd
/usr/local/share/cpufreqd/files/cpufreqd.conf
/usr/local/tce.installed/cpufreqd
