# tinycorelinux/16.x/x86/tcz/google-talkplugin.tcz.list
etc/cron.daily/google-talkplugin
opt/google/chrome/pepper/libppo1d.so
opt/google/talkplugin/GoogleTalkPlugin
opt/google/talkplugin/attributions.txt
opt/google/talkplugin/cron/google-talkplugin
opt/google/talkplugin/data/LMspeed_510.emd
opt/google/talkplugin/data/SFTprec_120.emd
opt/google/talkplugin/lib/libCg.so
opt/google/talkplugin/lib/libCgGL.so
opt/google/talkplugin/libgoogletalkremoting.so
opt/google/talkplugin/libnpgoogletalk.so
opt/google/talkplugin/libnpgtpo3dautoplugin.so
opt/google/talkplugin/libnpo1d.so
opt/google/talkplugin/libppo1d.so
opt/google/talkplugin/locale/af/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/am/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/ar/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/bg/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/bn/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/ca/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/cs/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/da/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/de/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/el/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/en/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/en-GB/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/es/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/es-419/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/et/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/eu/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/fa/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/fi/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/fil/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/fr/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/fr-CA/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/gl/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/gu/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/hi/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/hr/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/hu/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/id/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/is/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/it/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/iw/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/ja/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/kn/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/ko/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/lt/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/lv/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/ml/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/mr/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/ms/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/nl/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/no/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/or/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/pl/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/pt-BR/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/pt-PT/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/ro/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/ru/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/sk/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/sl/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/sr/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/sv/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/sw/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/ta/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/te/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/th/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/tl/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/tr/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/uk/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/ur/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/vi/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/zh-CN/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/zh-HK/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/zh-TW/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/locale/zu/LC_MESSAGES/windowpicker.mo
opt/google/talkplugin/remoting24x24.png
opt/google/talkplugin/windowpicker.glade
usr/lib/firefox/plugins/libnpgoogletalk.so
usr/lib/firefox/plugins/libnpgtpo3dautoplugin.so
usr/lib/firefox/plugins/libnpo1d.so
usr/lib/iceape/plugins/libnpgoogletalk.so
usr/lib/iceape/plugins/libnpgtpo3dautoplugin.so
usr/lib/iceape/plugins/libnpo1d.so
usr/lib/iceweasel/plugins/libnpgoogletalk.so
usr/lib/iceweasel/plugins/libnpgtpo3dautoplugin.so
usr/lib/iceweasel/plugins/libnpo1d.so
usr/lib/midbrowser/plugins/libnpgoogletalk.so
usr/lib/midbrowser/plugins/libnpgtpo3dautoplugin.so
usr/lib/midbrowser/plugins/libnpo1d.so
usr/lib/mozilla/plugins/libnpgoogletalk.so
usr/lib/mozilla/plugins/libnpgtpo3dautoplugin.so
usr/lib/mozilla/plugins/libnpo1d.so
usr/lib/xulrunner/plugins/libnpgoogletalk.so
usr/lib/xulrunner/plugins/libnpgtpo3dautoplugin.so
usr/lib/xulrunner/plugins/libnpo1d.so
usr/lib/xulrunner-addons/plugins/libnpgoogletalk.so
usr/lib/xulrunner-addons/plugins/libnpgtpo3dautoplugin.so
usr/lib/xulrunner-addons/plugins/libnpo1d.so
usr/local/tce.installed/google-talkplugin
