# tinycorelinux/16.x/x86/tcz/hostapd.tcz.list
/usr/local/tce.installed/hostapd
/usr/local/share/doc/hostapd/COPYING
/usr/local/etc/init.d/hostapd
/usr/local/bin/hostapd_cli
/usr/local/bin/hostapd
