# tinycorelinux/16.x/x86/tcz/no-ip.tcz.list
usr/local/bin/noip2
usr/local/etc/init.d/no-ip
usr/local/share/doc/License/no-ip.txt
