# tinycorelinux/16.x/x86/tcz/sc101-nbd.tcz.list
usr/local/etc/init.d/sc101-nbd
usr/local/etc/sc101-nbd.conf.sample
usr/local/sbin/ut
usr/local/share/man/man8/ut.8.gz
usr/local/tce.installed/sc101-nbd
