# tinycorelinux/16.x/x86/tcz/Xorg-7.7-3d-dev.tcz.list
/usr/local/include/GLES/egl.h
/usr/local/include/GLES/gl.h
/usr/local/include/GLES/glext.h
/usr/local/include/GLES/glplatform.h
/usr/local/include/GLES3/gl3.h
/usr/local/include/GLES3/gl31.h
/usr/local/include/GLES3/gl32.h
/usr/local/include/GLES3/gl3ext.h
/usr/local/include/GLES3/gl3platform.h
/usr/local/include/xa_composite.h
/usr/local/include/xa_context.h
/usr/local/include/xa_tracker.h
/usr/local/lib/pkgconfig/dri.pc
/usr/local/lib/pkgconfig/glesv1_cm.pc
/usr/local/lib/pkgconfig/xatracker.pc
/usr/local/share/drirc.d/00-mesa-defaults.conf
/usr/local/share/drirc.d/00-radv-defaults.conf
/usr/local/share/vulkan/icd.d/intel_hasvk_icd.i686.json
/usr/local/share/vulkan/icd.d/intel_icd.i686.json
/usr/local/share/vulkan/icd.d/lvp_icd.i686.json
/usr/local/share/vulkan/icd.d/radeon_icd.i686.json
