# tinycorelinux/13.x/aarch64/tcz/miniupnpd.tcz.list
/etc/init.d/miniupnpd
/etc/miniupnpd/ip6tables_init.sh
/etc/miniupnpd/ip6tables_removeall.sh
/etc/miniupnpd/iptables_init.sh
/etc/miniupnpd/iptables_removeall.sh
/etc/miniupnpd/miniupnpd.conf
/etc/miniupnpd/miniupnpd_functions.sh
/usr/local/sbin/miniupnpd
/usr/local/tce.installed/miniupnpd
