# tinycorelinux/8.x/x86_64/tcz/bind.tcz.list
/usr/local/bin/arpaname
/usr/local/bin/bind9-config
/usr/local/bin/delv
/usr/local/bin/dig
/usr/local/bin/host
/usr/local/bin/isc-config.sh
/usr/local/bin/mdig
/usr/local/bin/named-rrchecker
/usr/local/bin/nslookup
/usr/local/bin/nsupdate
/usr/local/etc/bind.keys
/usr/local/etc/init.d/bind
/usr/local/lib/libbind9.so
/usr/local/lib/libbind9.so.160
/usr/local/lib/libbind9.so.160.0.4
/usr/local/lib/libdns.so
/usr/local/lib/libdns.so.169
/usr/local/lib/libdns.so.169.0.1
/usr/local/lib/libirs.so
/usr/local/lib/libirs.so.160
/usr/local/lib/libirs.so.160.0.3
/usr/local/lib/libisccc.so
/usr/local/lib/libisccc.so.160
/usr/local/lib/libisccc.so.160.0.2
/usr/local/lib/libisccfg.so
/usr/local/lib/libisccfg.so.160
/usr/local/lib/libisccfg.so.160.1.1
/usr/local/lib/libisc.so
/usr/local/lib/libisc.so.166
/usr/local/lib/libisc.so.166.2.0
/usr/local/lib/liblwres.so
/usr/local/lib/liblwres.so.160
/usr/local/lib/liblwres.so.160.0.0
/usr/local/sbin/ddns-confgen
/usr/local/sbin/dnssec-dsfromkey
/usr/local/sbin/dnssec-importkey
/usr/local/sbin/dnssec-keyfromlabel
/usr/local/sbin/dnssec-keygen
/usr/local/sbin/dnssec-revoke
/usr/local/sbin/dnssec-settime
/usr/local/sbin/dnssec-signzone
/usr/local/sbin/dnssec-verify
/usr/local/sbin/genrandom
/usr/local/sbin/isc-hmac-fixup
/usr/local/sbin/lwresd
/usr/local/sbin/named
/usr/local/sbin/named-checkconf
/usr/local/sbin/named-checkzone
/usr/local/sbin/named-compilezone
/usr/local/sbin/named-journalprint
/usr/local/sbin/nsec3hash
/usr/local/sbin/rndc
/usr/local/sbin/rndc-confgen
/usr/local/sbin/tsig-keygen
