# tinycorelinux/8.x/x86_64/tcz/open-vm-tools.tcz.list
/etc/vmware-tools/resume-vm-default
/etc/vmware-tools/suspend-vm-default
/etc/vmware-tools/scripts/poweroff-vm-default.d/99-shutdown
/etc/vmware-tools/scripts/vmware/network
/etc/vmware-tools/poweron-vm-default
/etc/vmware-tools/statechange.subr
/etc/vmware-tools/guestproxy-ssl.conf
/etc/vmware-tools/poweroff-vm-default
/etc/vmware-tools/vm-support
/usr/local/etc/init.d/open-vm-tools
/usr/local/etc/udev/rules.d/99-vmware-scsi-udev.rules
/usr/local/sbin/mount.vmhgfs
/usr/local/lib/libguestlib.so
/usr/local/lib/libhgfs.so.0.0.0
/usr/local/lib/libhgfs.so.0
/usr/local/lib/libvmtools.so
/usr/local/lib/libguestlib.so.0
/usr/local/lib/libguestlib.so.0.0.0
/usr/local/lib/libvmtools.so.0.0.0
/usr/local/lib/open-vm-tools/plugins/common/libvix.so
/usr/local/lib/open-vm-tools/plugins/common/libhgfsServer.so
/usr/local/lib/open-vm-tools/plugins/vmsvc/libguestInfo.so
/usr/local/lib/open-vm-tools/plugins/vmsvc/libtimeSync.so
/usr/local/lib/open-vm-tools/plugins/vmsvc/libpowerOps.so
/usr/local/lib/open-vm-tools/plugins/vmsvc/libgrabbitmqProxy.so
/usr/local/lib/open-vm-tools/plugins/vmsvc/libresolutionKMS.so
/usr/local/lib/open-vm-tools/plugins/vmsvc/libvmbackup.so
/usr/local/lib/libvmtools.so.0
/usr/local/lib/libhgfs.so
/usr/local/tce.installed/open-vm-tools
/usr/local/bin/vmware-xferlogs
/usr/local/bin/vmhgfs-fuse
/usr/local/bin/vmware-rpctool
/usr/local/bin/vmware-hgfsclient
/usr/local/bin/vmware-guestproxycerttool
/usr/local/bin/vmware-vmblock-fuse
/usr/local/bin/vmware-toolbox-cmd
/usr/local/bin/vmtoolsd
/usr/local/bin/vmware-checkvm
/usr/local/bin/vmware-namespace-cmd
