# tinycorelinux/8.x/x86_64/tcz/samba.tcz.list
/usr/local/bin/cifsdd
/usr/local/bin/dbwrap_tool
/usr/local/bin/eventlogadm
/usr/local/bin/findsmb
/usr/local/bin/gentest
/usr/local/bin/ldbadd
/usr/local/bin/ldbdel
/usr/local/bin/ldbedit
/usr/local/bin/ldbmodify
/usr/local/bin/ldbrename
/usr/local/bin/ldbsearch
/usr/local/bin/locktest
/usr/local/bin/masktest
/usr/local/bin/mvxattr
/usr/local/bin/ndrdump
/usr/local/bin/net
/usr/local/bin/nmblookup
/usr/local/bin/ntlm_auth
/usr/local/bin/oLschema2ldif
/usr/local/bin/pdbedit
/usr/local/bin/pidl
/usr/local/bin/profiles
/usr/local/bin/regdiff
/usr/local/bin/regpatch
/usr/local/bin/regshell
/usr/local/bin/regtree
/usr/local/bin/rpcclient
/usr/local/bin/samba-regedit
/usr/local/bin/sharesec
/usr/local/bin/smbcacls
/usr/local/bin/smbclient
/usr/local/bin/smbcontrol
/usr/local/bin/smbcquotas
/usr/local/bin/smbget
/usr/local/bin/smbpasswd
/usr/local/bin/smbspool
/usr/local/bin/smbstatus
/usr/local/bin/smbtar
/usr/local/bin/smbtorture
/usr/local/bin/smbtree
/usr/local/bin/tdbbackup
/usr/local/bin/tdbdump
/usr/local/bin/tdbrestore
/usr/local/bin/tdbtool
/usr/local/bin/testparm
/usr/local/bin/wbinfo
/usr/local/etc/init.d/samba
/usr/local/etc/samba/smb.conf.default
/usr/local/lib/libdcerpc-binding.so
/usr/local/lib/libdcerpc-binding.so.0
/usr/local/lib/libdcerpc-binding.so.0.0.1
/usr/local/lib/libdcerpc-samr.so
/usr/local/lib/libdcerpc-samr.so.0
/usr/local/lib/libdcerpc-samr.so.0.0.1
/usr/local/lib/libdcerpc.so
/usr/local/lib/libdcerpc.so.0
/usr/local/lib/libdcerpc.so.0.0.1
/usr/local/lib/libndr-krb5pac.so
/usr/local/lib/libndr-krb5pac.so.0
/usr/local/lib/libndr-krb5pac.so.0.0.1
/usr/local/lib/libndr-nbt.so
/usr/local/lib/libndr-nbt.so.0
/usr/local/lib/libndr-nbt.so.0.0.1
/usr/local/lib/libndr-standard.so
/usr/local/lib/libndr-standard.so.0
/usr/local/lib/libndr-standard.so.0.0.1
/usr/local/lib/libndr.so
/usr/local/lib/libndr.so.0
/usr/local/lib/libndr.so.0.0.8
/usr/local/lib/libnetapi.so
/usr/local/lib/libnetapi.so.0
/usr/local/lib/libnss_winbind.so
/usr/local/lib/libnss_winbind.so.2
/usr/local/lib/libnss_wins.so
/usr/local/lib/libnss_wins.so.2
/usr/local/lib/libsamba-credentials.so
/usr/local/lib/libsamba-credentials.so.0
/usr/local/lib/libsamba-credentials.so.0.0.1
/usr/local/lib/libsamba-errors.so
/usr/local/lib/libsamba-errors.so.1
/usr/local/lib/libsamba-hostconfig.so
/usr/local/lib/libsamba-hostconfig.so.0
/usr/local/lib/libsamba-hostconfig.so.0.0.1
/usr/local/lib/libsamba-passdb.so
/usr/local/lib/libsamba-passdb.so.0
/usr/local/lib/libsamba-passdb.so.0.26.0
/usr/local/lib/libsamba-policy.so
/usr/local/lib/libsamba-policy.so.0
/usr/local/lib/libsamba-policy.so.0.0.1
/usr/local/lib/libsamba-util.so
/usr/local/lib/libsamba-util.so.0
/usr/local/lib/libsamba-util.so.0.0.1
/usr/local/lib/libsamdb.so
/usr/local/lib/libsamdb.so.0
/usr/local/lib/libsamdb.so.0.0.1
/usr/local/lib/libsmbclient.so
/usr/local/lib/libsmbclient.so.0
/usr/local/lib/libsmbclient.so.0.2.3
/usr/local/lib/libsmbconf.so
/usr/local/lib/libsmbconf.so.0
/usr/local/lib/libtevent-util.so
/usr/local/lib/libtevent-util.so.0
/usr/local/lib/libtevent-util.so.0.0.1
/usr/local/lib/samba/auth/script.so
/usr/local/lib/samba/idmap/autorid.so
/usr/local/lib/samba/idmap/hash.so
/usr/local/lib/samba/idmap/rid.so
/usr/local/lib/samba/idmap/script.so
/usr/local/lib/samba/idmap/tdb2.so
/usr/local/lib/samba/ldb/asq.so
/usr/local/lib/samba/ldb/ildap.so
/usr/local/lib/samba/ldb/ldbsamba_extensions.so
/usr/local/lib/samba/ldb/paged_results.so
/usr/local/lib/samba/ldb/paged_searches.so
/usr/local/lib/samba/ldb/rdn_name.so
/usr/local/lib/samba/ldb/sample.so
/usr/local/lib/samba/ldb/server_sort.so
/usr/local/lib/samba/ldb/skel.so
/usr/local/lib/samba/ldb/tdb.so
/usr/local/lib/samba/libCHARSET3-samba4.so
/usr/local/lib/samba/libLIBWBCLIENT-OLD-samba4.so
/usr/local/lib/samba/libMESSAGING-samba4.so
/usr/local/lib/samba/libaddns-samba4.so
/usr/local/lib/samba/libads-samba4.so
/usr/local/lib/samba/libasn1-samba4.so.8
/usr/local/lib/samba/libasn1-samba4.so.8.0.0
/usr/local/lib/samba/libasn1util-samba4.so
/usr/local/lib/samba/libauth-sam-reply-samba4.so
/usr/local/lib/samba/libauth-samba4.so
/usr/local/lib/samba/libauth-unix-token-samba4.so
/usr/local/lib/samba/libauth4-samba4.so
/usr/local/lib/samba/libauthkrb5-samba4.so
/usr/local/lib/samba/libcli-cldap-samba4.so
/usr/local/lib/samba/libcli-ldap-common-samba4.so
/usr/local/lib/samba/libcli-ldap-samba4.so
/usr/local/lib/samba/libcli-nbt-samba4.so
/usr/local/lib/samba/libcli-smb-common-samba4.so
/usr/local/lib/samba/libcli-spoolss-samba4.so
/usr/local/lib/samba/libcliauth-samba4.so
/usr/local/lib/samba/libcluster-samba4.so
/usr/local/lib/samba/libcmdline-credentials-samba4.so
/usr/local/lib/samba/libcom_err-samba4.so.0
/usr/local/lib/samba/libcom_err-samba4.so.0.25
/usr/local/lib/samba/libdbwrap-samba4.so
/usr/local/lib/samba/libdcerpc-samba-samba4.so
/usr/local/lib/samba/libdcerpc-samba4.so
/usr/local/lib/samba/libdsdb-garbage-collect-tombstones-samba4.so
/usr/local/lib/samba/libdsdb-module-samba4.so
/usr/local/lib/samba/libevents-samba4.so
/usr/local/lib/samba/libflag-mapping-samba4.so
/usr/local/lib/samba/libgenrand-samba4.so
/usr/local/lib/samba/libgensec-samba4.so
/usr/local/lib/samba/libgpo-samba4.so
/usr/local/lib/samba/libgse-samba4.so
/usr/local/lib/samba/libgssapi-samba4.so.2
/usr/local/lib/samba/libgssapi-samba4.so.2.0.0
/usr/local/lib/samba/libhcrypto-samba4.so.5
/usr/local/lib/samba/libhcrypto-samba4.so.5.0.1
/usr/local/lib/samba/libhdb-samba4.so.11
/usr/local/lib/samba/libhdb-samba4.so.11.0.2
/usr/local/lib/samba/libheimbase-samba4.so.1
/usr/local/lib/samba/libheimbase-samba4.so.1.0.0
/usr/local/lib/samba/libheimntlm-samba4.so.1
/usr/local/lib/samba/libheimntlm-samba4.so.1.0.1
/usr/local/lib/samba/libhttp-samba4.so
/usr/local/lib/samba/libhx509-samba4.so.5
/usr/local/lib/samba/libhx509-samba4.so.5.0.0
/usr/local/lib/samba/libidmap-samba4.so
/usr/local/lib/samba/libinterfaces-samba4.so
/usr/local/lib/samba/libiov-buf-samba4.so
/usr/local/lib/samba/libkdc-samba4.so.2
/usr/local/lib/samba/libkdc-samba4.so.2.0.0
/usr/local/lib/samba/libkrb5-samba4.so.26
/usr/local/lib/samba/libkrb5-samba4.so.26.0.0
/usr/local/lib/samba/libkrb5samba-samba4.so
/usr/local/lib/samba/libldb-cmdline-samba4.so
/usr/local/lib/samba/libldb.so.1
/usr/local/lib/samba/libldb.so.1.1.29
/usr/local/lib/samba/libldbsamba-samba4.so
/usr/local/lib/samba/liblibcli-lsa3-samba4.so
/usr/local/lib/samba/liblibcli-netlogon3-samba4.so
/usr/local/lib/samba/liblibsmb-samba4.so
/usr/local/lib/samba/libmessages-dgm-samba4.so
/usr/local/lib/samba/libmessages-util-samba4.so
/usr/local/lib/samba/libmsghdr-samba4.so
/usr/local/lib/samba/libmsrpc3-samba4.so
/usr/local/lib/samba/libndr-samba-samba4.so
/usr/local/lib/samba/libndr-samba4.so
/usr/local/lib/samba/libnet-keytab-samba4.so
/usr/local/lib/samba/libnetif-samba4.so
/usr/local/lib/samba/libnon-posix-acls-samba4.so
/usr/local/lib/samba/libnpa-tstream-samba4.so
/usr/local/lib/samba/libnss-info-samba4.so
/usr/local/lib/samba/libpopt-samba3-samba4.so
/usr/local/lib/samba/libpopt-samba4.so
/usr/local/lib/samba/libposix-eadb-samba4.so
/usr/local/lib/samba/libprinting-migrate-samba4.so
/usr/local/lib/samba/libpyldb-util.so.1
/usr/local/lib/samba/libpyldb-util.so.1.1.29
/usr/local/lib/samba/libpytalloc-util.so.2
/usr/local/lib/samba/libpytalloc-util.so.2.1.9
/usr/local/lib/samba/libregistry-samba4.so
/usr/local/lib/samba/libroken-samba4.so.19
/usr/local/lib/samba/libroken-samba4.so.19.0.1
/usr/local/lib/samba/libsamba-cluster-support-samba4.so
/usr/local/lib/samba/libsamba-debug-samba4.so
/usr/local/lib/samba/libsamba-modules-samba4.so
/usr/local/lib/samba/libsamba-net-samba4.so
/usr/local/lib/samba/libsamba-python-samba4.so
/usr/local/lib/samba/libsamba-security-samba4.so
/usr/local/lib/samba/libsamba-sockets-samba4.so
/usr/local/lib/samba/libsamba3-util-samba4.so
/usr/local/lib/samba/libsamdb-common-samba4.so
/usr/local/lib/samba/libsecrets3-samba4.so
/usr/local/lib/samba/libserver-id-db-samba4.so
/usr/local/lib/samba/libserver-role-samba4.so
/usr/local/lib/samba/libshares-samba4.so
/usr/local/lib/samba/libsmb-transport-samba4.so
/usr/local/lib/samba/libsmbclient-raw-samba4.so
/usr/local/lib/samba/libsmbd-base-samba4.so
/usr/local/lib/samba/libsmbd-conn-samba4.so
/usr/local/lib/samba/libsmbd-shim-samba4.so
/usr/local/lib/samba/libsmbpasswdparser-samba4.so
/usr/local/lib/samba/libsocket-blocking-samba4.so
/usr/local/lib/samba/libsys-rw-samba4.so
/usr/local/lib/samba/libtalloc-report-samba4.so
/usr/local/lib/samba/libtalloc.so.2
/usr/local/lib/samba/libtalloc.so.2.1.9
/usr/local/lib/samba/libtdb-wrap-samba4.so
/usr/local/lib/samba/libtdb.so.1
/usr/local/lib/samba/libtdb.so.1.3.12
/usr/local/lib/samba/libtevent.so.0
/usr/local/lib/samba/libtevent.so.0.9.31
/usr/local/lib/samba/libtime-basic-samba4.so
/usr/local/lib/samba/libtorture-samba4.so
/usr/local/lib/samba/libtrusts-util-samba4.so
/usr/local/lib/samba/libutil-cmdline-samba4.so
/usr/local/lib/samba/libutil-reg-samba4.so
/usr/local/lib/samba/libutil-setid-samba4.so
/usr/local/lib/samba/libutil-tdb-samba4.so
/usr/local/lib/samba/libwind-samba4.so.0
/usr/local/lib/samba/libwind-samba4.so.0.0.0
/usr/local/lib/samba/libxattr-tdb-samba4.so
/usr/local/lib/samba/nss_info/hash.so
/usr/local/lib/samba/samba/smbspool_krb5_wrapper
/usr/local/lib/samba/vfs/acl_tdb.so
/usr/local/lib/samba/vfs/acl_xattr.so
/usr/local/lib/samba/vfs/aio_fork.so
/usr/local/lib/samba/vfs/aio_pthread.so
/usr/local/lib/samba/vfs/audit.so
/usr/local/lib/samba/vfs/btrfs.so
/usr/local/lib/samba/vfs/cap.so
/usr/local/lib/samba/vfs/catia.so
/usr/local/lib/samba/vfs/commit.so
/usr/local/lib/samba/vfs/crossrename.so
/usr/local/lib/samba/vfs/default_quota.so
/usr/local/lib/samba/vfs/dirsort.so
/usr/local/lib/samba/vfs/expand_msdfs.so
/usr/local/lib/samba/vfs/extd_audit.so
/usr/local/lib/samba/vfs/fake_acls.so
/usr/local/lib/samba/vfs/fake_perms.so
/usr/local/lib/samba/vfs/fileid.so
/usr/local/lib/samba/vfs/fruit.so
/usr/local/lib/samba/vfs/full_audit.so
/usr/local/lib/samba/vfs/linux_xfs_sgid.so
/usr/local/lib/samba/vfs/media_harmony.so
/usr/local/lib/samba/vfs/netatalk.so
/usr/local/lib/samba/vfs/nfs4acl_xattr.so
/usr/local/lib/samba/vfs/offline.so
/usr/local/lib/samba/vfs/preopen.so
/usr/local/lib/samba/vfs/readahead.so
/usr/local/lib/samba/vfs/readonly.so
/usr/local/lib/samba/vfs/recycle.so
/usr/local/lib/samba/vfs/shadow_copy.so
/usr/local/lib/samba/vfs/shadow_copy2.so
/usr/local/lib/samba/vfs/shell_snap.so
/usr/local/lib/samba/vfs/snapper.so
/usr/local/lib/samba/vfs/streams_depot.so
/usr/local/lib/samba/vfs/streams_xattr.so
/usr/local/lib/samba/vfs/syncops.so
/usr/local/lib/samba/vfs/time_audit.so
/usr/local/lib/samba/vfs/unityed_media.so
/usr/local/lib/samba/vfs/worm.so
/usr/local/lib/samba/vfs/xattr_tdb.so
/usr/local/lib/winbind_krb5_locator.so
/usr/local/sbin/nmbd
/usr/local/sbin/smbd
/usr/local/sbin/winbindd
/usr/local/share/samba/files/smb.conf
/usr/local/tce.installed/samba
