# tinycorelinux/16.x/aarch64/tcz/bluez.tcz.list
lib/udev/hid2hci
usr/local/tce.installed/bluez
usr/local/share/dbus-1/services/org.bluez.obex.service
usr/local/share/bluez/files/97-hid2hci.rules
usr/local/share/bluez/files/proximity.conf
usr/local/share/bluez/files/network.conf
usr/local/share/bluez/files/main.conf
usr/local/share/bluez/files/input.conf
usr/local/lib/cups/backend/bluetooth
usr/local/lib/bluez/test/test-sap-server
usr/local/lib/bluez/test/test-profile
usr/local/lib/bluez/test/test-network
usr/local/lib/bluez/test/test-nap
usr/local/lib/bluez/test/test-mesh
usr/local/lib/bluez/test/test-manager
usr/local/lib/bluez/test/test-hfp
usr/local/lib/bluez/test/test-health-sink
usr/local/lib/bluez/test/test-health
usr/local/lib/bluez/test/test-gatt-profile
usr/local/lib/bluez/test/test-discovery
usr/local/lib/bluez/test/test-device
usr/local/lib/bluez/test/test-adapter
usr/local/lib/bluez/test/simple-player
usr/local/lib/bluez/test/simple-endpoint
usr/local/lib/bluez/test/simple-agent
usr/local/lib/bluez/test/service-spp.xml
usr/local/lib/bluez/test/service-record.dtd
usr/local/lib/bluez/test/service-opp.xml
usr/local/lib/bluez/test/service-ftp.xml
usr/local/lib/bluez/test/service-did.xml
usr/local/lib/bluez/test/sap_client.py
usr/local/lib/bluez/test/pbap-client
usr/local/lib/bluez/test/opp-client
usr/local/lib/bluez/test/monitor-bluetooth
usr/local/lib/bluez/test/map-client
usr/local/lib/bluez/test/list-devices
usr/local/lib/bluez/test/ftp-client
usr/local/lib/bluez/test/example-gatt-server
usr/local/lib/bluez/test/example-gatt-client
usr/local/lib/bluez/test/example-advertisement
usr/local/lib/bluez/test/dbusdef.py
usr/local/lib/bluez/test/bluezutils.py
usr/local/lib/bluez/test/agent.py
usr/local/lib/bluetooth/obexd
usr/local/lib/bluetooth/bluetoothd
usr/local/lib/bluetooth/plugins/sixaxis.so
usr/local/etc/init.d/bluez
usr/local/etc/dbus-1/system.d/bluetooth.conf
usr/local/bin/sdptool
usr/local/bin/rfcomm
usr/local/bin/rctest
usr/local/bin/mpris-proxy
usr/local/bin/l2test
usr/local/bin/l2ping
usr/local/bin/isotest
usr/local/bin/hex2hcd
usr/local/bin/hcitool
usr/local/bin/hcidump
usr/local/bin/hciconfig
usr/local/bin/hciattach
usr/local/bin/gatttool
usr/local/bin/ciptool
usr/local/bin/btmon
usr/local/bin/btattach
usr/local/bin/bluetoothctl
usr/local/bin/bluemoon
