# tinycorelinux/16.x/aarch64/tcz/mesa.tcz.list
/usr/local/lib/dri/kms_swrast_dri.so
/usr/local/lib/dri/libdril_dri.so
/usr/local/lib/dri/swrast_dri.so
/usr/local/lib/dri/v3d_dri.so
/usr/local/lib/dri/vc4_dri.so
/usr/local/lib/dri/virtio_gpu_dri.so
/usr/local/lib/dri/vkms_dri.so
/usr/local/lib/libGL.so
/usr/local/lib/libGL.so.1
/usr/local/lib/libGL.so.1.2.0
/usr/local/lib/libgallium-24.2.5.so
/usr/local/lib/libglapi.so
/usr/local/lib/libglapi.so.0
/usr/local/lib/libglapi.so.0.0.0
/usr/local/share/drirc.d/00-mesa-defaults.conf
/usr/local/tce.installed/mesa
