# tinycorelinux/16.x/aarch64/tcz/prometheus.tcz.list
/usr/local/etc/init.d/prometheus
/usr/local/etc/prometheus/prometheus_example.yml
/usr/local/bin/promtool
/usr/local/bin/prometheus
/usr/local/share/doc/prometheus/LICENSE
/usr/local/share/doc/prometheus/NOTICE
/usr/local/tce.installed/prometheus
