# tinycorelinux/12.x/x86/tcz/frr.tcz.list
/usr/local/tce.installed/frr
/usr/local/sbin/zebra
/usr/local/sbin/watchfrr.sh
/usr/local/sbin/watchfrr
/usr/local/sbin/staticd
/usr/local/sbin/ssd
/usr/local/sbin/ripngd
/usr/local/sbin/ripd
/usr/local/sbin/pimd
/usr/local/sbin/pbrd
/usr/local/sbin/ospfd
/usr/local/sbin/ospf6d
/usr/local/sbin/nhrpd
/usr/local/sbin/ldpd
/usr/local/sbin/isisd
/usr/local/sbin/frrinit.sh
/usr/local/sbin/frrcommon.sh
/usr/local/sbin/frr-reload.py
/usr/local/sbin/frr-reload
/usr/local/sbin/eigrpd
/usr/local/sbin/bgpd
/usr/local/sbin/bfdd
/usr/local/sbin/babeld
/usr/local/lib/libfrrospfapiclient.so.0.0.0
/usr/local/lib/libfrrospfapiclient.so.0
/usr/local/lib/libfrrospfapiclient.so
/usr/local/lib/libfrr.so.0.0.0
/usr/local/lib/libfrr.so.0
/usr/local/lib/libfrr.so
/usr/local/lib/frr/modules/zebra_irdp.so
/usr/local/etc/init.d/frr
/usr/local/etc/frr/zebra.conf.sample
/usr/local/etc/frr/vtysh.conf.sample
/usr/local/etc/frr/staticd.conf.sample
/usr/local/etc/frr/ripngd.conf.sample
/usr/local/etc/frr/ripd.conf.sample
/usr/local/etc/frr/pimd.conf.sample
/usr/local/etc/frr/pbrd.conf.sample
/usr/local/etc/frr/ospfd.conf.sample
/usr/local/etc/frr/ospf6d.conf.sample
/usr/local/etc/frr/ldpd.conf.sample
/usr/local/etc/frr/isisd.conf.sample
/usr/local/etc/frr/eigrpd.conf.sample
/usr/local/etc/frr/daemons.sample
/usr/local/etc/frr/bgpd.conf.vnc.sample
/usr/local/etc/frr/bgpd.conf.sample2
/usr/local/etc/frr/bgpd.conf.sample
/usr/local/etc/frr/bfdd.conf.sample
/usr/local/etc/frr/babeld.conf.sample
/usr/local/bin/vtysh
/usr/local/bin/mtracebis
