# tinycorelinux/12.x/x86/tcz/pcsc-lite.tcz.list
/usr/local/etc/init.d/pcscd
/usr/local/lib/libpcsclite.so
/usr/local/lib/libpcsclite.so.1
/usr/local/lib/libpcsclite.so.1.0.0
/usr/local/sbin/pcscd
