# tinycorelinux/12.x/x86/tcz/pdnsd.tcz.list
usr/local/etc/init.d/pdnsd
usr/local/etc/pdnsd.conf.default
usr/local/etc/pdnsd.conf.sample
usr/local/sbin/pdnsd
usr/local/sbin/pdnsd-ctl
usr/local/tce.installed/pdnsd
