# tinycorelinux/6.x/x86/tcz/bluez.tcz.list
/lib/udev/bluetooth_serial
/usr/local/bin/ciptool
/usr/local/bin/dund
/usr/local/bin/hcitool
/usr/local/bin/hidd
/usr/local/bin/l2ping
/usr/local/bin/pand
/usr/local/bin/rfcomm
/usr/local/bin/sdptool
/usr/local/bin/simple-agent
/usr/local/bin/simple-service
/usr/local/bin/test-discovery
/usr/local/bin/test-network
/usr/local/etc/alsa/example_.asoundrc
/usr/local/etc/dbus-1/system.d/bluetooth.conf
/usr/local/etc/init.d/bluez
/usr/local/lib/alsa-lib/libasound_module_ctl_bluetooth.so
/usr/local/lib/alsa-lib/libasound_module_pcm_bluetooth.so
/usr/local/lib/cups/backend/bluetooth
/usr/local/lib/gstreamer-0.10/libgstbluetooth.so
/usr/local/lib/libbluetooth.so
/usr/local/lib/libbluetooth.so.3
/usr/local/lib/libbluetooth.so.3.13.0
/usr/local/sbin/bccmd
/usr/local/sbin/bluetoothd
/usr/local/sbin/hciattach
/usr/local/sbin/hciconfig
/usr/local/share/alsa/bluetooth.conf
/usr/local/share/bluez/files/97-bluetooth-serial.rules
/usr/local/share/bluez/files/audio.conf
/usr/local/share/bluez/files/input.conf
/usr/local/share/bluez/files/main.conf
/usr/local/share/bluez/files/network.conf
/usr/local/share/bluez/files/proximity.conf
/usr/local/share/bluez/files/rfcomm.conf
/usr/local/share/bluez/files/serial.conf
/usr/local/share/dbus-1/system-services/org.bluez.service
/usr/local/tce.installed/bluez
