# tinycorelinux/6.x/x86/tcz/lighttpd.tcz.list
usr/local/etc/init.d/lighttp
usr/local/lib/lighttpd/mod_access.so
usr/local/lib/lighttpd/mod_accesslog.so
usr/local/lib/lighttpd/mod_alias.so
usr/local/lib/lighttpd/mod_auth.so
usr/local/lib/lighttpd/mod_cgi.so
usr/local/lib/lighttpd/mod_cml.so
usr/local/lib/lighttpd/mod_compress.so
usr/local/lib/lighttpd/mod_dirlisting.so
usr/local/lib/lighttpd/mod_evasive.so
usr/local/lib/lighttpd/mod_evhost.so
usr/local/lib/lighttpd/mod_expire.so
usr/local/lib/lighttpd/mod_extforward.so
usr/local/lib/lighttpd/mod_fastcgi.so
usr/local/lib/lighttpd/mod_flv_streaming.so
usr/local/lib/lighttpd/mod_indexfile.so
usr/local/lib/lighttpd/mod_magnet.so
usr/local/lib/lighttpd/mod_mysql_vhost.so
usr/local/lib/lighttpd/mod_proxy.so
usr/local/lib/lighttpd/mod_redirect.so
usr/local/lib/lighttpd/mod_rewrite.so
usr/local/lib/lighttpd/mod_rrdtool.so
usr/local/lib/lighttpd/mod_scgi.so
usr/local/lib/lighttpd/mod_secdownload.so
usr/local/lib/lighttpd/mod_setenv.so
usr/local/lib/lighttpd/mod_simple_vhost.so
usr/local/lib/lighttpd/mod_ssi.so
usr/local/lib/lighttpd/mod_staticfile.so
usr/local/lib/lighttpd/mod_status.so
usr/local/lib/lighttpd/mod_trigger_b4_dl.so
usr/local/lib/lighttpd/mod_userdir.so
usr/local/lib/lighttpd/mod_usertrack.so
usr/local/lib/lighttpd/mod_webdav.so
usr/local/sbin/lighttpd
usr/local/sbin/lighttpd-angel
usr/local/share/doc/License/lighttp-copying
usr/local/share/doc/lighttpd/conf.d/access_log.conf
usr/local/share/doc/lighttpd/conf.d/auth.conf
usr/local/share/doc/lighttpd/conf.d/cgi.conf
usr/local/share/doc/lighttpd/conf.d/cml.conf
usr/local/share/doc/lighttpd/conf.d/compress.conf
usr/local/share/doc/lighttpd/conf.d/debug.conf
usr/local/share/doc/lighttpd/conf.d/dirlisting.conf
usr/local/share/doc/lighttpd/conf.d/evhost.conf
usr/local/share/doc/lighttpd/conf.d/expire.conf
usr/local/share/doc/lighttpd/conf.d/fastcgi.conf
usr/local/share/doc/lighttpd/conf.d/geoip.conf
usr/local/share/doc/lighttpd/conf.d/magnet.conf
usr/local/share/doc/lighttpd/conf.d/mime.conf
usr/local/share/doc/lighttpd/conf.d/mod.template
usr/local/share/doc/lighttpd/conf.d/mysql_vhost.conf
usr/local/share/doc/lighttpd/conf.d/proxy.conf
usr/local/share/doc/lighttpd/conf.d/rrdtool.conf
usr/local/share/doc/lighttpd/conf.d/scgi.conf
usr/local/share/doc/lighttpd/conf.d/secdownload.conf
usr/local/share/doc/lighttpd/conf.d/simple_vhost.conf
usr/local/share/doc/lighttpd/conf.d/ssi.conf
usr/local/share/doc/lighttpd/conf.d/status.conf
usr/local/share/doc/lighttpd/conf.d/trigger_b4_dl.conf
usr/local/share/doc/lighttpd/conf.d/userdir.conf
usr/local/share/doc/lighttpd/conf.d/webdav.conf
usr/local/share/doc/lighttpd/lighttpd.conf
usr/local/share/doc/lighttpd/modules.conf
usr/local/share/doc/lighttpd/vhosts.d/template.conf
usr/local/tce.installed/lighttpd
