# tinycorelinux/6.x/x86/tcz/open-vm-tools-3.16.6-tinycore.tcz.list
etc/init.d/open-vm-tools
etc/open-vm-tools/guestproxy-ssl.conf
etc/open-vm-tools/locations
etc/open-vm-tools/plugins
etc/open-vm-tools/poweroff-vm-default
etc/open-vm-tools/poweron-vm-default
etc/open-vm-tools/resume-vm-default
etc/open-vm-tools/scripts/open-vm-tools/network
etc/open-vm-tools/statechange.subr
etc/open-vm-tools/suspend-vm-default
etc/open-vm-tools/vm-support
etc/release
lib/modules/3.16.6-tinycore/kernel/fs/vmhgfs/vmhgfs.ko
sbin/mount.vmhgfs
usr/local/bin/vmtoolsd
usr/local/bin/vmware-checkvm
usr/local/bin/vmware-guestproxycerttool
usr/local/bin/vmware-hgfsclient
usr/local/bin/vmware-rpctool
usr/local/bin/vmware-toolbox-cmd
usr/local/bin/vmware-user
usr/local/bin/vmware-user-suid-wrapper
usr/local/bin/vmware-xferlogs
usr/local/etc/xdg/autostart/vmware-user.desktop
usr/local/include/libDeployPkg/guestcust-events.h
usr/local/include/libDeployPkg/imgcust-api.h
usr/local/include/libDeployPkg/includeCheck.h
usr/local/include/libDeployPkg/linuxDeployment.h
usr/local/include/libDeployPkg/log.h
usr/local/include/libDeployPkg/process.h
usr/local/include/libDeployPkg/rpcout.h
usr/local/include/libDeployPkg/vm_basic_types.h
usr/local/include/vmGuestLib/includeCheck.h
usr/local/include/vmGuestLib/vmGuestLib.h
usr/local/include/vmGuestLib/vmSessionId.h
usr/local/include/vmGuestLib/vm_basic_types.h
usr/local/include/vmci_sockets.h
usr/local/lib/libDeployPkg.la
usr/local/lib/libDeployPkg.so
usr/local/lib/libDeployPkg.so.0
usr/local/lib/libDeployPkg.so.0.0.0
usr/local/lib/libguestlib.la
usr/local/lib/libguestlib.so
usr/local/lib/libguestlib.so.0
usr/local/lib/libguestlib.so.0.0.0
usr/local/lib/libhgfs.la
usr/local/lib/libhgfs.so
usr/local/lib/libhgfs.so.0
usr/local/lib/libhgfs.so.0.0.0
usr/local/lib/libvmtools.la
usr/local/lib/libvmtools.so
usr/local/lib/libvmtools.so.0
usr/local/lib/libvmtools.so.0.0.0
usr/local/lib/open-vm-tools/bin/vmtoolsd
usr/local/lib/open-vm-tools/bin/vmware-checkvm
usr/local/lib/open-vm-tools/bin/vmware-guestproxycerttool
usr/local/lib/open-vm-tools/bin/vmware-hgfsclient
usr/local/lib/open-vm-tools/bin/vmware-rpctool
usr/local/lib/open-vm-tools/bin/vmware-toolbox-cmd
usr/local/lib/open-vm-tools/bin/vmware-user-suid-wrapper
usr/local/lib/open-vm-tools/bin/vmware-xferlogs
usr/local/lib/open-vm-tools/lib/libDeployPkg.la
usr/local/lib/open-vm-tools/lib/libDeployPkg.so
usr/local/lib/open-vm-tools/lib/libDeployPkg.so.0
usr/local/lib/open-vm-tools/lib/libDeployPkg.so.0.0.0
usr/local/lib/open-vm-tools/lib/libguestlib.la
usr/local/lib/open-vm-tools/lib/libguestlib.so
usr/local/lib/open-vm-tools/lib/libguestlib.so.0
usr/local/lib/open-vm-tools/lib/libguestlib.so.0.0.0
usr/local/lib/open-vm-tools/lib/libhgfs.la
usr/local/lib/open-vm-tools/lib/libhgfs.so
usr/local/lib/open-vm-tools/lib/libhgfs.so.0
usr/local/lib/open-vm-tools/lib/libhgfs.so.0.0.0
usr/local/lib/open-vm-tools/lib/libvmtools.la
usr/local/lib/open-vm-tools/lib/libvmtools.so
usr/local/lib/open-vm-tools/lib/libvmtools.so.0
usr/local/lib/open-vm-tools/lib/libvmtools.so.0.0.0
usr/local/lib/open-vm-tools/plugins/common/libhgfsServer.la
usr/local/lib/open-vm-tools/plugins/common/libhgfsServer.so
usr/local/lib/open-vm-tools/plugins/common/libvix.la
usr/local/lib/open-vm-tools/plugins/common/libvix.so
usr/local/lib/open-vm-tools/plugins/vmsvc/libdeployPkgPlugin.so
usr/local/lib/open-vm-tools/plugins/vmsvc/libgrabbitmqProxy.so
usr/local/lib/open-vm-tools/plugins/vmsvc/libguestInfo.so
usr/local/lib/open-vm-tools/plugins/vmsvc/libpowerOps.so
usr/local/lib/open-vm-tools/plugins/vmsvc/libtimeSync.so
usr/local/lib/open-vm-tools/plugins/vmsvc/libvmbackup.so
usr/local/lib/open-vm-tools/plugins/vmusr/libdesktopEvents.so
usr/local/lib/open-vm-tools/plugins/vmusr/libdndcp.so
usr/local/lib/open-vm-tools/plugins/vmusr/libresolutionSet.so
usr/local/lib/open-vm-tools/sbin/mount.vmhgfs
usr/local/lib/pkgconfig/libDeployPkg.pc
usr/local/lib/pkgconfig/pkgconfig
usr/local/lib/pkgconfig/vmguestlib.pc
usr/local/sbin/mount.vmhgfs
usr/local/share/open-vm-tools/messages/de/toolboxcmd.vmsg
usr/local/share/open-vm-tools/messages/de/vmtoolsd.vmsg
usr/local/share/open-vm-tools/messages/ja/toolboxcmd.vmsg
usr/local/share/open-vm-tools/messages/ja/vmtoolsd.vmsg
usr/local/share/open-vm-tools/messages/ko/toolboxcmd.vmsg
usr/local/share/open-vm-tools/messages/ko/vmtoolsd.vmsg
usr/local/share/open-vm-tools/messages/zh_CN/toolboxcmd.vmsg
usr/local/tce.installed/open-vm-tools-3.16.6-tinycore
