# tinycorelinux/6.x/x86/tcz/smartmontools.tcz.list
/usr/local/etc/init.d/smartd
/usr/local/etc/smartd.conf.sample
/usr/local/etc/smartd_warning.sh
/usr/local/sbin/smartctl
/usr/local/sbin/smartd
/usr/local/sbin/update-smart-drivedb
/usr/local/share/smartmontools/files/drivedb.h
/usr/local/share/smartmontools/files/smartd.conf
/usr/local/tce.installed/smartmontools
