# tinycorelinux/7.x/x86_64/tcz/acpid.tcz.list
/usr/local/bin/acpi_listen
/usr/local/etc/init.d/acpid
/usr/local/sbin/acpid
/usr/local/sbin/kacpimon
/usr/local/share/acpid/files/lid
/usr/local/share/acpid/files/lid.sh
/usr/local/tce.installed/acpid
