# tinycorelinux/7.x/x86_64/tcz/alsa-config.tcz.list
/usr/local/bin/speaker-test
/usr/local/etc/alsa/alsa.conf.d/README
/usr/local/etc/alsa/cards/AACI.conf
/usr/local/etc/alsa/cards/ATIIXP-MODEM.conf
/usr/local/etc/alsa/cards/ATIIXP-SPDMA.conf
/usr/local/etc/alsa/cards/ATIIXP.conf
/usr/local/etc/alsa/cards/AU8810.conf
/usr/local/etc/alsa/cards/AU8820.conf
/usr/local/etc/alsa/cards/AU8830.conf
/usr/local/etc/alsa/cards/Audigy.conf
/usr/local/etc/alsa/cards/Audigy2.conf
/usr/local/etc/alsa/cards/Aureon51.conf
/usr/local/etc/alsa/cards/Aureon71.conf
/usr/local/etc/alsa/cards/CA0106.conf
/usr/local/etc/alsa/cards/CMI8338-SWIEC.conf
/usr/local/etc/alsa/cards/CMI8338.conf
/usr/local/etc/alsa/cards/CMI8738-MC6.conf
/usr/local/etc/alsa/cards/CMI8738-MC8.conf
/usr/local/etc/alsa/cards/CMI8788.conf
/usr/local/etc/alsa/cards/CS46xx.conf
/usr/local/etc/alsa/cards/EMU10K1.conf
/usr/local/etc/alsa/cards/EMU10K1X.conf
/usr/local/etc/alsa/cards/ENS1370.conf
/usr/local/etc/alsa/cards/ENS1371.conf
/usr/local/etc/alsa/cards/ES1968.conf
/usr/local/etc/alsa/cards/Echo_Echo3G.conf
/usr/local/etc/alsa/cards/FM801.conf
/usr/local/etc/alsa/cards/FWSpeakers.conf
/usr/local/etc/alsa/cards/FireWave.conf
/usr/local/etc/alsa/cards/GUS.conf
/usr/local/etc/alsa/cards/HDA-Intel.conf
/usr/local/etc/alsa/cards/ICE1712.conf
/usr/local/etc/alsa/cards/ICE1724.conf
/usr/local/etc/alsa/cards/ICH-MODEM.conf
/usr/local/etc/alsa/cards/ICH.conf
/usr/local/etc/alsa/cards/ICH4.conf
/usr/local/etc/alsa/cards/Loopback.conf
/usr/local/etc/alsa/cards/Maestro3.conf
/usr/local/etc/alsa/cards/NFORCE.conf
/usr/local/etc/alsa/cards/PC-Speaker.conf
/usr/local/etc/alsa/cards/PMac.conf
/usr/local/etc/alsa/cards/PMacToonie.conf
/usr/local/etc/alsa/cards/PS3.conf
/usr/local/etc/alsa/cards/RME9636.conf
/usr/local/etc/alsa/cards/RME9652.conf
/usr/local/etc/alsa/cards/SB-XFi.conf
/usr/local/etc/alsa/cards/SI7018/sndoc-mixer.alisp
/usr/local/etc/alsa/cards/SI7018/sndop-mixer.alisp
/usr/local/etc/alsa/cards/SI7018.conf
/usr/local/etc/alsa/cards/TRID4DWAVENX.conf
/usr/local/etc/alsa/cards/USB-Audio.conf
/usr/local/etc/alsa/cards/VIA686A.conf
/usr/local/etc/alsa/cards/VIA8233.conf
/usr/local/etc/alsa/cards/VIA8233A.conf
/usr/local/etc/alsa/cards/VIA8237.conf
/usr/local/etc/alsa/cards/VX222.conf
/usr/local/etc/alsa/cards/VXPocket.conf
/usr/local/etc/alsa/cards/VXPocket440.conf
/usr/local/etc/alsa/cards/YMF744.conf
/usr/local/etc/alsa/cards/aliases.alisp
/usr/local/etc/alsa/cards/aliases.conf
/usr/local/etc/alsa/pcm/center_lfe.conf
/usr/local/etc/alsa/pcm/default.conf
/usr/local/etc/alsa/pcm/dmix.conf
/usr/local/etc/alsa/pcm/dpl.conf
/usr/local/etc/alsa/pcm/dsnoop.conf
/usr/local/etc/alsa/pcm/front.conf
/usr/local/etc/alsa/pcm/hdmi.conf
/usr/local/etc/alsa/pcm/iec958.conf
/usr/local/etc/alsa/pcm/modem.conf
/usr/local/etc/alsa/pcm/rear.conf
/usr/local/etc/alsa/pcm/side.conf
/usr/local/etc/alsa/pcm/surround21.conf
/usr/local/etc/alsa/pcm/surround40.conf
/usr/local/etc/alsa/pcm/surround41.conf
/usr/local/etc/alsa/pcm/surround50.conf
/usr/local/etc/alsa/pcm/surround51.conf
/usr/local/etc/alsa/pcm/surround71.conf
/usr/local/etc/alsa/topology/broadwell/broadwell.conf
/usr/local/etc/alsa/ucm/DAISY-I2S/DAISY-I2S.conf
/usr/local/etc/alsa/ucm/DAISY-I2S/HiFi.conf
/usr/local/etc/alsa/ucm/GoogleNyan
/usr/local/etc/alsa/ucm/GoogleNyan/GoogleNyan.conf
/usr/local/etc/alsa/ucm/GoogleNyan/HiFi.conf
/usr/local/etc/alsa/ucm/PAZ00/HiFi.conf
/usr/local/etc/alsa/ucm/PAZ00/PAZ00.conf
/usr/local/etc/alsa/ucm/PAZ00/Record.conf
/usr/local/etc/alsa/ucm/PandaBoard/FMAnalog
/usr/local/etc/alsa/ucm/PandaBoard/PandaBoard.conf
/usr/local/etc/alsa/ucm/PandaBoard/hifi
/usr/local/etc/alsa/ucm/PandaBoard/hifiLP
/usr/local/etc/alsa/ucm/PandaBoard/record
/usr/local/etc/alsa/ucm/PandaBoard/voice
/usr/local/etc/alsa/ucm/PandaBoard/voiceCall
/usr/local/etc/alsa/ucm/PandaBoardES/FMAnalog
/usr/local/etc/alsa/ucm/PandaBoardES/PandaBoardES.conf
/usr/local/etc/alsa/ucm/PandaBoardES/hifi
/usr/local/etc/alsa/ucm/PandaBoardES/hifiLP
/usr/local/etc/alsa/ucm/PandaBoardES/record
/usr/local/etc/alsa/ucm/PandaBoardES/voice
/usr/local/etc/alsa/ucm/PandaBoardES/voiceCall
/usr/local/etc/alsa/ucm/SDP4430/FMAnalog
/usr/local/etc/alsa/ucm/SDP4430/SDP4430.conf
/usr/local/etc/alsa/ucm/SDP4430/hifi
/usr/local/etc/alsa/ucm/SDP4430/hifiLP
/usr/local/etc/alsa/ucm/SDP4430/record
/usr/local/etc/alsa/ucm/SDP4430/voice
/usr/local/etc/alsa/ucm/SDP4430/voiceCall
/usr/local/etc/alsa/ucm/broadwell-rt286
/usr/local/etc/alsa/ucm/broadwell-rt286/HiFi
/usr/local/etc/alsa/ucm/broadwell-rt286/broadwell-rt286.conf
/usr/local/etc/alsa/ucm/tegraalc5632/tegraalc5632.conf
/usr/local/etc/init.d/alsasound
/usr/local/sbin/alsa-info.sh
/usr/local/share/alsa/init/00main
/usr/local/share/alsa/init/ca0106
/usr/local/share/alsa/init/default
/usr/local/share/alsa/init/hda
/usr/local/share/alsa/init/help
/usr/local/share/alsa/init/info
/usr/local/share/alsa/init/test
/usr/local/share/alsa/speaker-test/sample_map.csv
/usr/local/share/sounds/alsa/Front_Center.wav
/usr/local/share/sounds/alsa/Front_Left.wav
/usr/local/share/sounds/alsa/Front_Right.wav
/usr/local/share/sounds/alsa/Noise.wav
/usr/local/share/sounds/alsa/Rear_Center.wav
/usr/local/share/sounds/alsa/Rear_Left.wav
/usr/local/share/sounds/alsa/Rear_Right.wav
/usr/local/share/sounds/alsa/Side_Left.wav
/usr/local/share/sounds/alsa/Side_Right.wav
