# tinycorelinux/7.x/x86_64/tcz/audit.tcz.list
/usr/local/bin/aulast
/usr/local/bin/aulastlog
/usr/local/bin/ausyscall
/usr/local/bin/auvirt
/usr/local/etc/audisp/audisp-remote.conf
/usr/local/etc/audisp/audispd.conf
/usr/local/etc/audisp/plugins.d/af_unix.conf
/usr/local/etc/audisp/plugins.d/au-remote.conf
/usr/local/etc/audisp/plugins.d/audispd-zos-remote.conf
/usr/local/etc/audisp/plugins.d/syslog.conf
/usr/local/etc/audisp/zos-remote.conf
/usr/local/etc/audit/audit-stop.rules
/usr/local/etc/audit/auditd.conf
/usr/local/etc/libaudit.conf
/usr/local/etc/rc.d/init.d/auditd
/usr/local/etc/sysconfig/auditd
/usr/local/lib/libaudit.so
/usr/local/lib/libaudit.so.1
/usr/local/lib/libaudit.so.1.0.0
/usr/local/lib/libauparse.so
/usr/local/lib/libauparse.so.0
/usr/local/lib/libauparse.so.0.0.0
/usr/local/lib/python2.7/site-packages/_audit.so
/usr/local/lib/python2.7/site-packages/audit.py
/usr/local/lib/python2.7/site-packages/auparse.so
/usr/local/sbin/audisp-remote
/usr/local/sbin/audispd
/usr/local/sbin/audispd-zos-remote
/usr/local/sbin/auditctl
/usr/local/sbin/auditd
/usr/local/sbin/augenrules
/usr/local/sbin/aureport
/usr/local/sbin/ausearch
/usr/local/sbin/autrace
