# tinycorelinux/7.x/x86_64/tcz/pcsc-lite.tcz.list
/usr/local/bin/pcsc-spy
/usr/local/etc/init.d/pcscd
/usr/local/lib/libpcsclite.so
/usr/local/lib/libpcsclite.so.1
/usr/local/lib/libpcsclite.so.1.0.0
/usr/local/lib/libpcscspy.so
/usr/local/lib/libpcscspy.so.0
/usr/local/lib/libpcscspy.so.0.0.0
/usr/local/sbin/pcscd
