# tinycorelinux/7.x/armv6/tcz/dropbear.tcz.list
usr/local/etc/init.d/dropbear
usr/local/etc/dropbear/banner
usr/local/sbin/dropbear
usr/local/bin/dropbearmulti
usr/local/bin/scp
usr/local/bin/dropbearconvert
usr/local/bin/dropbearkey
usr/local/bin/dbclient
usr/local/bin/ssh
