# tinycorelinux/12.x/x86_64/tcz/bind.tcz.list
/usr/local/bin/arpaname
/usr/local/bin/delv
/usr/local/bin/dig
/usr/local/bin/host
/usr/local/bin/mdig
/usr/local/bin/named-rrchecker
/usr/local/bin/nslookup
/usr/local/bin/nsupdate
/usr/local/etc/bind.keys
/usr/local/etc/init.d/bind
/usr/local/lib/libbind9.so
/usr/local/lib/libbind9.so.1600
/usr/local/lib/libbind9.so.1600.0.6
/usr/local/lib/libdns.so
/usr/local/lib/libdns.so.1610
/usr/local/lib/libdns.so.1610.0.0
/usr/local/lib/libirs.so
/usr/local/lib/libirs.so.1601
/usr/local/lib/libirs.so.1601.0.0
/usr/local/lib/libisccc.so
/usr/local/lib/libisccc.so.1600
/usr/local/lib/libisccc.so.1600.0.2
/usr/local/lib/libisccfg.so
/usr/local/lib/libisccfg.so.1602
/usr/local/lib/libisccfg.so.1602.0.0
/usr/local/lib/libisc.so
/usr/local/lib/libisc.so.1608
/usr/local/lib/libisc.so.1608.0.1
/usr/local/lib/libns.so
/usr/local/lib/libns.so.1606
/usr/local/lib/libns.so.1606.0.1
/usr/local/lib/named/filter-aaaa.so
/usr/local/lib/python3.6/site-packages/isc-2.0-py3.6.egg-info
/usr/local/lib/python3.6/site-packages/isc/checkds.py
/usr/local/lib/python3.6/site-packages/isc/coverage.py
/usr/local/lib/python3.6/site-packages/isc/dnskey.py
/usr/local/lib/python3.6/site-packages/isc/eventlist.py
/usr/local/lib/python3.6/site-packages/isc/__init__.py
/usr/local/lib/python3.6/site-packages/isc/keydict.py
/usr/local/lib/python3.6/site-packages/isc/keyevent.py
/usr/local/lib/python3.6/site-packages/isc/keymgr.py
/usr/local/lib/python3.6/site-packages/isc/keyseries.py
/usr/local/lib/python3.6/site-packages/isc/keyzone.py
/usr/local/lib/python3.6/site-packages/isc/parsetab.py
/usr/local/lib/python3.6/site-packages/isc/policy.py
/usr/local/lib/python3.6/site-packages/isc/__pycache__/checkds.cpython-36.pyc
/usr/local/lib/python3.6/site-packages/isc/__pycache__/coverage.cpython-36.pyc
/usr/local/lib/python3.6/site-packages/isc/__pycache__/dnskey.cpython-36.pyc
/usr/local/lib/python3.6/site-packages/isc/__pycache__/eventlist.cpython-36.pyc
/usr/local/lib/python3.6/site-packages/isc/__pycache__/__init__.cpython-36.pyc
/usr/local/lib/python3.6/site-packages/isc/__pycache__/keydict.cpython-36.pyc
/usr/local/lib/python3.6/site-packages/isc/__pycache__/keyevent.cpython-36.pyc
/usr/local/lib/python3.6/site-packages/isc/__pycache__/keymgr.cpython-36.pyc
/usr/local/lib/python3.6/site-packages/isc/__pycache__/keyseries.cpython-36.pyc
/usr/local/lib/python3.6/site-packages/isc/__pycache__/keyzone.cpython-36.pyc
/usr/local/lib/python3.6/site-packages/isc/__pycache__/parsetab.cpython-36.pyc
/usr/local/lib/python3.6/site-packages/isc/__pycache__/policy.cpython-36.pyc
/usr/local/lib/python3.6/site-packages/isc/__pycache__/rndc.cpython-36.pyc
/usr/local/lib/python3.6/site-packages/isc/__pycache__/utils.cpython-36.pyc
/usr/local/lib/python3.6/site-packages/isc/rndc.py
/usr/local/lib/python3.6/site-packages/isc/utils.py
/usr/local/sbin/ddns-confgen
/usr/local/sbin/dnssec-cds
/usr/local/sbin/dnssec-checkds
/usr/local/sbin/dnssec-coverage
/usr/local/sbin/dnssec-dsfromkey
/usr/local/sbin/dnssec-importkey
/usr/local/sbin/dnssec-keyfromlabel
/usr/local/sbin/dnssec-keygen
/usr/local/sbin/dnssec-keymgr
/usr/local/sbin/dnssec-revoke
/usr/local/sbin/dnssec-settime
/usr/local/sbin/dnssec-signzone
/usr/local/sbin/dnssec-verify
/usr/local/sbin/named
/usr/local/sbin/named-checkconf
/usr/local/sbin/named-checkzone
/usr/local/sbin/named-compilezone
/usr/local/sbin/named-journalprint
/usr/local/sbin/nsec3hash
/usr/local/sbin/rndc
/usr/local/sbin/rndc-confgen
/usr/local/sbin/tsig-keygen
