# tinycorelinux/12.x/x86_64/tcz/hostapd.tcz.list
/usr/local/etc/init.d/hostapd
/usr/local/tce.installed/hostapd
/usr/local/bin/hostapd
/usr/local/bin/hostapd_cli
