# tinycorelinux/14.x/x86_64/tcz/mongo.tcz.list
/usr/local/bin/mongo
/usr/local/etc/init.d/mongodb
/usr/local/etc/mongod/mongod.conf.default
/usr/local/sbin/mongod
/usr/local/sbin/mongos
/usr/local/share/doc/mongo-5.3.0/LICENSE-Community.txt
/usr/local/share/doc/mongo-5.3.0/MPL-2
/usr/local/share/doc/mongo-5.3.0/THIRD-PARTY-NOTICES
/usr/local/share/mongod/files/mongod.conf
/usr/local/tce.installed/mongo
