# tinycorelinux/6.x/armv7/tcz/dbus.tcz.list
usr/local/tce.installed/dbus
usr/local/etc/init.d/dbus
usr/local/etc/dbus-1/system.conf
usr/local/etc/dbus-1/session.conf
usr/local/libexec/dbus-daemon-launch-helper
usr/local/bin/dbus-daemon
usr/local/bin/dbus-launch
usr/local/bin/dbus-monitor
usr/local/bin/dbus-send
usr/local/bin/dbus-cleanup-sockets
usr/local/bin/dbus-run-session
usr/local/bin/dbus-uuidgen
usr/local/lib/libdbus-1.so.3.8.11
usr/local/lib/libdbus-1.so
usr/local/lib/libdbus-1.so.3
