# tinycorelinux/11.x/x86/tcz/at.tcz.list
/usr/local/bin/at
/usr/local/bin/atq
/usr/local/bin/atrm
/usr/local/bin/batch
/usr/local/etc/init.d/atd
/usr/local/sbin/atd
/usr/local/sbin/atrun
/usr/local/share/at/files/at.deny
/usr/local/share/doc/at/Copyright
/usr/local/tce.installed/at
