# tinycorelinux/11.x/x86/tcz/p910nd.tcz.list
/usr/local/tce.installed/p910nd
/usr/local/sbin/p910nd
/usr/local/etc/p910nd.conf.sample
/usr/local/etc/init.d/p910nd
