# tinycorelinux/5.x/x86_64/tcz/cups.tcz.list
/usr/local/bin/cancel
/usr/local/bin/cupstestdsc
/usr/local/bin/cupstestppd
/usr/local/bin/ipptool
/usr/local/bin/lp
/usr/local/bin/lpoptions
/usr/local/bin/lppasswd
/usr/local/bin/lpq
/usr/local/bin/lpr
/usr/local/bin/lprm
/usr/local/bin/lpstat
/usr/local/bin/ppdc
/usr/local/bin/ppdhtml
/usr/local/bin/ppdi
/usr/local/bin/ppdmerge
/usr/local/bin/ppdpo
/usr/local/etc/cups/cups-files.conf
/usr/local/etc/cups/cupsd.conf.default
/usr/local/etc/cups/snmp.conf
/usr/local/etc/init.d/cups
/usr/local/lib/cups/backend/dnssd
/usr/local/lib/cups/backend/http
/usr/local/lib/cups/backend/https
/usr/local/lib/cups/backend/ipp
/usr/local/lib/cups/backend/ipps
/usr/local/lib/cups/backend/lpd
/usr/local/lib/cups/backend/smb
/usr/local/lib/cups/backend/snmp
/usr/local/lib/cups/backend/socket
/usr/local/lib/cups/backend/usb
/usr/local/lib/cups/cgi-bin/admin.cgi
/usr/local/lib/cups/cgi-bin/classes.cgi
/usr/local/lib/cups/cgi-bin/help.cgi
/usr/local/lib/cups/cgi-bin/jobs.cgi
/usr/local/lib/cups/cgi-bin/printers.cgi
/usr/local/lib/cups/daemon/cups-deviced
/usr/local/lib/cups/daemon/cups-driverd
/usr/local/lib/cups/daemon/cups-exec
/usr/local/lib/cups/daemon/cups-lpd
/usr/local/lib/cups/driver
/usr/local/lib/cups/filter/commandtops
/usr/local/lib/cups/filter/gziptoany
/usr/local/lib/cups/filter/pstops
/usr/local/lib/cups/filter/rastertodymo
/usr/local/lib/cups/filter/rastertoepson
/usr/local/lib/cups/filter/rastertohp
/usr/local/lib/cups/filter/rastertolabel
/usr/local/lib/cups/filter/rastertopwg
/usr/local/lib/cups/monitor/bcp
/usr/local/lib/cups/monitor/tbcp
/usr/local/lib/cups/notifier/dbus
/usr/local/lib/cups/notifier/mailto
/usr/local/lib/cups/notifier/rss
/usr/local/sbin/accept
/usr/local/sbin/cupsaccept
/usr/local/sbin/cupsaddsmb
/usr/local/sbin/cupsctl
/usr/local/sbin/cupsd
/usr/local/sbin/cupsdisable
/usr/local/sbin/cupsenable
/usr/local/sbin/cupsfilter
/usr/local/sbin/cupsreject
/usr/local/sbin/lpadmin
/usr/local/sbin/lpc
/usr/local/sbin/lpinfo
/usr/local/sbin/lpmove
/usr/local/sbin/reject
/usr/local/share/cups/drv/sample.drv
/usr/local/share/cups/examples/color.drv
/usr/local/share/cups/examples/constraint.drv
/usr/local/share/cups/examples/custom.drv
/usr/local/share/cups/examples/grouping.drv
/usr/local/share/cups/examples/laserjet-basic.drv
/usr/local/share/cups/examples/laserjet-pjl.drv
/usr/local/share/cups/examples/minimum.drv
/usr/local/share/cups/examples/postscript.drv
/usr/local/share/cups/examples/r300-basic.drv
/usr/local/share/cups/examples/r300-colorman.drv
/usr/local/share/cups/examples/r300-remote.drv
/usr/local/share/cups/files/cupsd.conf
/usr/local/share/cups/ipptool/color.jpg
/usr/local/share/cups/ipptool/create-printer-subscription.test
/usr/local/share/cups/ipptool/document-a4.pdf
/usr/local/share/cups/ipptool/document-a4.ps
/usr/local/share/cups/ipptool/document-letter.pdf
/usr/local/share/cups/ipptool/document-letter.ps
/usr/local/share/cups/ipptool/get-completed-jobs.test
/usr/local/share/cups/ipptool/get-jobs.test
/usr/local/share/cups/ipptool/get-printer-attributes.test
/usr/local/share/cups/ipptool/gray.jpg
/usr/local/share/cups/ipptool/ipp-1.1.test
/usr/local/share/cups/ipptool/ipp-2.0.test
/usr/local/share/cups/ipptool/ipp-2.1.test
/usr/local/share/cups/ipptool/ipp-2.2.test
/usr/local/share/cups/ipptool/onepage-a4.pdf
/usr/local/share/cups/ipptool/onepage-a4.ps
/usr/local/share/cups/ipptool/onepage-letter.pdf
/usr/local/share/cups/ipptool/onepage-letter.ps
/usr/local/share/cups/ipptool/testfile.jpg
/usr/local/share/cups/ipptool/testfile.pcl
/usr/local/share/cups/ipptool/testfile.pdf
/usr/local/share/cups/ipptool/testfile.ps
/usr/local/share/cups/ipptool/testfile.txt
/usr/local/share/cups/mime/mime.convs
/usr/local/share/cups/mime/mime.types
/usr/local/share/cups/ppdc/epson.h
/usr/local/share/cups/ppdc/font.defs
/usr/local/share/cups/ppdc/hp.h
/usr/local/share/cups/ppdc/label.h
/usr/local/share/cups/ppdc/media.defs
/usr/local/share/cups/ppdc/raster.defs
/usr/local/share/cups/templates/add-class.tmpl
/usr/local/share/cups/templates/add-printer.tmpl
/usr/local/share/cups/templates/add-rss-subscription.tmpl
/usr/local/share/cups/templates/admin.tmpl
/usr/local/share/cups/templates/choose-device.tmpl
/usr/local/share/cups/templates/choose-make.tmpl
/usr/local/share/cups/templates/choose-model.tmpl
/usr/local/share/cups/templates/choose-serial.tmpl
/usr/local/share/cups/templates/choose-uri.tmpl
/usr/local/share/cups/templates/class-added.tmpl
/usr/local/share/cups/templates/class-confirm.tmpl
/usr/local/share/cups/templates/class-deleted.tmpl
/usr/local/share/cups/templates/class-jobs-header.tmpl
/usr/local/share/cups/templates/class-modified.tmpl
/usr/local/share/cups/templates/class.tmpl
/usr/local/share/cups/templates/classes-header.tmpl
/usr/local/share/cups/templates/classes.tmpl
/usr/local/share/cups/templates/command.tmpl
/usr/local/share/cups/templates/edit-config.tmpl
/usr/local/share/cups/templates/error-op.tmpl
/usr/local/share/cups/templates/error.tmpl
/usr/local/share/cups/templates/header.tmpl
/usr/local/share/cups/templates/help-header.tmpl
/usr/local/share/cups/templates/help-printable.tmpl
/usr/local/share/cups/templates/help-trailer.tmpl
/usr/local/share/cups/templates/job-cancel.tmpl
/usr/local/share/cups/templates/job-hold.tmpl
/usr/local/share/cups/templates/job-move.tmpl
/usr/local/share/cups/templates/job-moved.tmpl
/usr/local/share/cups/templates/job-release.tmpl
/usr/local/share/cups/templates/job-restart.tmpl
/usr/local/share/cups/templates/jobs-header.tmpl
/usr/local/share/cups/templates/jobs.tmpl
/usr/local/share/cups/templates/list-available-printers.tmpl
/usr/local/share/cups/templates/modify-class.tmpl
/usr/local/share/cups/templates/modify-printer.tmpl
/usr/local/share/cups/templates/norestart.tmpl
/usr/local/share/cups/templates/option-boolean.tmpl
/usr/local/share/cups/templates/option-conflict.tmpl
/usr/local/share/cups/templates/option-header.tmpl
/usr/local/share/cups/templates/option-pickmany.tmpl
/usr/local/share/cups/templates/option-pickone.tmpl
/usr/local/share/cups/templates/option-trailer.tmpl
/usr/local/share/cups/templates/pager.tmpl
/usr/local/share/cups/templates/printer-accept.tmpl
/usr/local/share/cups/templates/printer-added.tmpl
/usr/local/share/cups/templates/printer-configured.tmpl
/usr/local/share/cups/templates/printer-confirm.tmpl
/usr/local/share/cups/templates/printer-default.tmpl
/usr/local/share/cups/templates/printer-deleted.tmpl
/usr/local/share/cups/templates/printer-jobs-header.tmpl
/usr/local/share/cups/templates/printer-modified.tmpl
/usr/local/share/cups/templates/printer-purge.tmpl
/usr/local/share/cups/templates/printer-reject.tmpl
/usr/local/share/cups/templates/printer-start.tmpl
/usr/local/share/cups/templates/printer-stop.tmpl
/usr/local/share/cups/templates/printer.tmpl
/usr/local/share/cups/templates/printers-header.tmpl
/usr/local/share/cups/templates/printers.tmpl
/usr/local/share/cups/templates/restart.tmpl
/usr/local/share/cups/templates/samba-export.tmpl
/usr/local/share/cups/templates/samba-exported.tmpl
/usr/local/share/cups/templates/search.tmpl
/usr/local/share/cups/templates/set-printer-options-header.tmpl
/usr/local/share/cups/templates/set-printer-options-trailer.tmpl
/usr/local/share/cups/templates/subscription-added.tmpl
/usr/local/share/cups/templates/subscription-canceled.tmpl
/usr/local/share/cups/templates/test-page.tmpl
/usr/local/share/cups/templates/trailer.tmpl
/usr/local/share/cups/templates/users.tmpl
/usr/local/share/doc/cups/cups-printable.css
/usr/local/share/doc/cups/cups.css
/usr/local/share/doc/cups/help/accounting.html
/usr/local/share/doc/cups/help/api-array.html
/usr/local/share/doc/cups/help/api-cgi.html
/usr/local/share/doc/cups/help/api-cups.html
/usr/local/share/doc/cups/help/api-driver.html
/usr/local/share/doc/cups/help/api-filedir.html
/usr/local/share/doc/cups/help/api-filter.html
/usr/local/share/doc/cups/help/api-httpipp.html
/usr/local/share/doc/cups/help/api-mime.html
/usr/local/share/doc/cups/help/api-overview.html
/usr/local/share/doc/cups/help/api-ppd.html
/usr/local/share/doc/cups/help/api-ppdc.html
/usr/local/share/doc/cups/help/api-raster.html
/usr/local/share/doc/cups/help/cgi.html
/usr/local/share/doc/cups/help/glossary.html
/usr/local/share/doc/cups/help/kerberos.html
/usr/local/share/doc/cups/help/license.html
/usr/local/share/doc/cups/help/man-backend.html
/usr/local/share/doc/cups/help/man-cancel.html
/usr/local/share/doc/cups/help/man-cups-config.html
/usr/local/share/doc/cups/help/man-cups-lpd.html
/usr/local/share/doc/cups/help/man-cups-snmp.html
/usr/local/share/doc/cups/help/man-cupsaccept.html
/usr/local/share/doc/cups/help/man-cupsaddsmb.html
/usr/local/share/doc/cups/help/man-cupsd.html
/usr/local/share/doc/cups/help/man-cupsenable.html
/usr/local/share/doc/cups/help/man-cupstestdsc.html
/usr/local/share/doc/cups/help/man-cupstestppd.html
/usr/local/share/doc/cups/help/man-filter.html
/usr/local/share/doc/cups/help/man-ipptool.html
/usr/local/share/doc/cups/help/man-ipptoolfile.html
/usr/local/share/doc/cups/help/man-lp.html
/usr/local/share/doc/cups/help/man-lpadmin.html
/usr/local/share/doc/cups/help/man-lpc.html
/usr/local/share/doc/cups/help/man-lpinfo.html
/usr/local/share/doc/cups/help/man-lpmove.html
/usr/local/share/doc/cups/help/man-lpoptions.html
/usr/local/share/doc/cups/help/man-lppasswd.html
/usr/local/share/doc/cups/help/man-lpq.html
/usr/local/share/doc/cups/help/man-lpr.html
/usr/local/share/doc/cups/help/man-lprm.html
/usr/local/share/doc/cups/help/man-lpstat.html
/usr/local/share/doc/cups/help/man-mime.convs.html
/usr/local/share/doc/cups/help/man-mime.types.html
/usr/local/share/doc/cups/help/man-notifier.html
/usr/local/share/doc/cups/help/man-ppdc.html
/usr/local/share/doc/cups/help/man-ppdhtml.html
/usr/local/share/doc/cups/help/man-ppdi.html
/usr/local/share/doc/cups/help/man-ppdmerge.html
/usr/local/share/doc/cups/help/man-ppdpo.html
/usr/local/share/doc/cups/help/network.html
/usr/local/share/doc/cups/help/options.html
/usr/local/share/doc/cups/help/overview.html
/usr/local/share/doc/cups/help/policies.html
/usr/local/share/doc/cups/help/postscript-driver.html
/usr/local/share/doc/cups/help/ppd-compiler.html
/usr/local/share/doc/cups/help/raster-driver.html
/usr/local/share/doc/cups/help/ref-access_log.html
/usr/local/share/doc/cups/help/ref-classes-conf.html
/usr/local/share/doc/cups/help/ref-client-conf.html
/usr/local/share/doc/cups/help/ref-cups-files-conf.html
/usr/local/share/doc/cups/help/ref-cupsd-conf.html
/usr/local/share/doc/cups/help/ref-error_log.html
/usr/local/share/doc/cups/help/ref-mailto-conf.html
/usr/local/share/doc/cups/help/ref-page_log.html
/usr/local/share/doc/cups/help/ref-ppdcfile.html
/usr/local/share/doc/cups/help/ref-printers-conf.html
/usr/local/share/doc/cups/help/ref-snmp-conf.html
/usr/local/share/doc/cups/help/ref-subscriptions-conf.html
/usr/local/share/doc/cups/help/security.html
/usr/local/share/doc/cups/help/sharing.html
/usr/local/share/doc/cups/help/spec-banner.html
/usr/local/share/doc/cups/help/spec-browsing.html
/usr/local/share/doc/cups/help/spec-cmp.html
/usr/local/share/doc/cups/help/spec-command.html
/usr/local/share/doc/cups/help/spec-design.html
/usr/local/share/doc/cups/help/spec-ipp.html
/usr/local/share/doc/cups/help/spec-pdf.html
/usr/local/share/doc/cups/help/spec-postscript.html
/usr/local/share/doc/cups/help/spec-ppd.html
/usr/local/share/doc/cups/help/spec-raster.html
/usr/local/share/doc/cups/help/spec-stp.html
/usr/local/share/doc/cups/help/standard.html
/usr/local/share/doc/cups/help/translation.html
/usr/local/share/doc/cups/help/whatsnew.html
/usr/local/share/doc/cups/images/color-wheel.png
/usr/local/share/doc/cups/images/cups-block-diagram.png
/usr/local/share/doc/cups/images/cups-command-chain.png
/usr/local/share/doc/cups/images/cups-icon.png
/usr/local/share/doc/cups/images/cups-postscript-chain.png
/usr/local/share/doc/cups/images/cups-raster-chain.png
/usr/local/share/doc/cups/images/cups.png
/usr/local/share/doc/cups/images/left.gif
/usr/local/share/doc/cups/images/raster-organization.png
/usr/local/share/doc/cups/images/raster.png
/usr/local/share/doc/cups/images/right.gif
/usr/local/share/doc/cups/images/sample-image.png
/usr/local/share/doc/cups/images/sel.gif
/usr/local/share/doc/cups/images/smiley.jpg
/usr/local/share/doc/cups/images/unsel.gif
/usr/local/share/doc/cups/images/wait.gif
/usr/local/share/doc/cups/index.html
/usr/local/share/doc/cups/robots.txt
/usr/local/tce.installed/cups
