# tinycorelinux/5.x/x86_64/tcz/networkmanager.tcz.list
/usr/local/bin/nm-online
/usr/local/bin/nm-tool
/usr/local/bin/nmcli
/usr/local/etc/NetworkManager/nm-system-settings.conf
/usr/local/etc/dbus-1/system.d/nm-avahi-autoipd.conf
/usr/local/etc/dbus-1/system.d/nm-dhcp-client.conf
/usr/local/etc/dbus-1/system.d/nm-dispatcher.conf
/usr/local/etc/dbus-1/system.d/org.freedesktop.NetworkManager.conf
/usr/local/etc/default/NetworkManager
/usr/local/etc/init.d/networkmanager
/usr/local/etc/polkit-1/rules.d/50-org.freedesktop.NetworkManager.rules
/usr/local/lib/NetworkManager/nm-avahi-autoipd.action
/usr/local/lib/NetworkManager/nm-dhcp-client.action
/usr/local/lib/NetworkManager/nm-dispatcher.action
/usr/local/lib/libnm-glib-vpn.so
/usr/local/lib/libnm-glib-vpn.so.1
/usr/local/lib/libnm-glib-vpn.so.1.1.0
/usr/local/lib/libnm-glib.so
/usr/local/lib/libnm-glib.so.4
/usr/local/lib/libnm-glib.so.4.6.0
/usr/local/lib/libnm-util.so
/usr/local/lib/libnm-util.so.2
/usr/local/lib/libnm-util.so.2.4.0
/usr/local/sbin/NetworkManager
/usr/local/share/bash-completion/completions/nmcli
/usr/local/share/dbus-1/system-services/org.freedesktop.nm_dispatcher.service
/usr/local/share/polkit-1/actions/org.freedesktop.NetworkManager.policy
/usr/local/tce.installed/networkmanager
