# tinycorelinux/5.x/armv6/tcz/dbus.tcz.list
usr/local/share/dbus/files/session.conf
usr/local/share/dbus/files/system.conf
usr/local/bin/dbus-send
usr/local/bin/dbus-run-session
usr/local/bin/dbus-cleanup-sockets
usr/local/bin/dbus-daemon
usr/local/bin/dbus-launch
usr/local/bin/dbus-uuidgen
usr/local/bin/dbus-monitor
usr/local/tce.installed/dbus
usr/local/libexec/dbus-daemon-launch-helper
usr/local/etc/init.d/dbus
usr/local/lib/libdbus-1.so
usr/local/lib/libdbus-1.so.3.8.3
usr/local/lib/libdbus-1.so.3
