# tinycorelinux/5.x/armv6/tcz/nfs-utils.tcz.list
usr/local/sbin/rpc.mountd
usr/local/sbin/nfsiostat
usr/local/sbin/osd_login
usr/local/sbin/mount.nfs
usr/local/sbin/start-statd
usr/local/sbin/rpc.statd
usr/local/sbin/umount.nfs
usr/local/sbin/sm-notify
usr/local/sbin/nfsstat
usr/local/sbin/mountstats
usr/local/sbin/showmount
usr/local/sbin/exportfs
usr/local/sbin/rpc.nfsd
usr/local/sbin/rpcdebug
usr/local/share/nfs-utils/files/exports
usr/local/tce.installed/nfs-utils
usr/local/etc/init.d/nfs-client
usr/local/etc/init.d/nfs-server
