# tinycorelinux/5.x/armv6/tcz/nginx.tcz.list
usr/local/tce.installed/nginx
usr/local/share/licences/nginx/LICENSE
usr/local/etc/init.d/nginx
usr/local/etc/nginx/nginx.conf.default
usr/local/etc/nginx/scgi_params.default
usr/local/etc/nginx/uwsgi_params.default
usr/local/etc/nginx/fastcgi.conf.default
usr/local/etc/nginx/fastcgi_params.default
usr/local/etc/nginx/mime.types.default
usr/local/etc/nginx/win-utf
usr/local/etc/nginx/koi-utf
usr/local/etc/nginx/koi-win
usr/local/sbin/nginx
usr/local/nginx/html/50x.html
usr/local/nginx/html/index.html
