# tinycorelinux/9.x/x86_64/tcz/bluez.tcz.list
/lib/udev/hid2hci
/usr/local/bin/bccmd
/usr/local/bin/bluemoon
/usr/local/bin/bluetoothctl
/usr/local/bin/btattach
/usr/local/bin/btmon
/usr/local/bin/hex2hcd
/usr/local/bin/l2ping
/usr/local/bin/l2test
/usr/local/bin/mpris-proxy
/usr/local/bin/rctest
/usr/local/etc/dbus-1/system.d/bluetooth.conf
/usr/local/etc/init.d/bluez
/usr/local/lib/bluetooth/bluetoothd
/usr/local/lib/bluetooth/obexd
/usr/local/lib/cups/backend/bluetooth
/usr/local/share/bluez/files/97-hid2hci.rules
/usr/local/share/bluez/files/input.conf
/usr/local/share/bluez/files/main.conf
/usr/local/share/bluez/files/network.conf
/usr/local/share/bluez/files/proximity.conf
/usr/local/share/bluez/files/rfcomm.conf
/usr/local/tce.installed/bluez
