# tinycorelinux/14.x/x86/tcz/PubSysFigures.tcz.list
/usr/local/tce.installed/PubSysFigures
/usr/local/sbin/PubSysFigures
/usr/local/etc/init.d/PubSybFigures
