# tinycorelinux/15.x/x86/tcz/minidlna.tcz.list
/usr/local/share/doc/minidlna/LICENCE.miniupnpd
/usr/local/share/doc/minidlna/COPYING
/usr/local/tce.installed/minidlna
/usr/local/sbin/minidlnad
/usr/local/etc/init.d/minidlna
/usr/local/etc/minidlna.conf
