# tinycorelinux/15.x/x86/tcz/proftpd.tcz.list
/usr/local/tce.installed/proftpd
/usr/local/share/man/man8/proftpd.8
/usr/local/share/man/man8/ftpshut.8
/usr/local/share/man/man8/ftpscrub.8
/usr/local/share/man/man8/ftpdctl.8
/usr/local/share/man/man5/xferlog.5
/usr/local/share/man/man5/proftpd.conf.5
/usr/local/share/man/man1/ftpwho.1
/usr/local/share/man/man1/ftptop.1
/usr/local/share/man/man1/ftpquota.1
/usr/local/share/man/man1/ftpmail.1
/usr/local/share/man/man1/ftpcount.1
/usr/local/share/man/man1/ftpasswd.1
/usr/local/share/doc/proftpd/COPYING
/usr/local/sbin/proftpd
/usr/local/sbin/ftpshut
/usr/local/sbin/ftpscrub
/usr/local/sbin/in.proftpd
/usr/local/lib/pkgconfig/proftpd.pc
/usr/local/etc/init.d/proftpd
/usr/local/bin/ftpwho
/usr/local/bin/ftptop
/usr/local/bin/ftpdctl
/usr/local/bin/ftpcount
/usr/local/bin/prxs
/usr/local/bin/ftpquota
/usr/local/bin/ftpmail
/usr/local/bin/ftpasswd
