# tinycorelinux/15.x/x86/tcz/pure-ftpd.tcz.list
/usr/local/tce.installed/pure-ftpd
/usr/local/share/man/man8/pure-uploadscript.8
/usr/local/share/man/man8/pure-statsdecode.8
/usr/local/share/man/man8/pure-quotacheck.8
/usr/local/share/man/man8/pure-pwconvert.8
/usr/local/share/man/man8/pure-pw.8
/usr/local/share/man/man8/pure-mrtginfo.8
/usr/local/share/man/man8/pure-ftpwho.8
/usr/local/share/man/man8/pure-ftpd.8
/usr/local/share/man/man8/pure-certd.8
/usr/local/share/man/man8/pure-authd.8
/usr/local/share/doc/pure-ftpd/COPYING
/usr/local/sbin/pure-uploadscript
/usr/local/sbin/pure-quotacheck
/usr/local/sbin/pure-mrtginfo
/usr/local/sbin/pure-ftpwho
/usr/local/sbin/pure-ftpd
/usr/local/sbin/pure-certd
/usr/local/sbin/pure-authd
/usr/local/etc/init.d/pure-ftpd
/usr/local/bin/pure-statsdecode
/usr/local/bin/pure-pwconvert
/usr/local/bin/pure-pw
