# tinycorelinux/10.x/x86/tcz/php7-fpm.tcz.list
/usr/local/etc/httpd/original/conf.d/httpd-php7-fpm.conf
/usr/local/etc/init.d/php-fpm
/usr/local/etc/php-fpm.conf.default
/usr/local/etc/php-fpm.d/www.conf.default
/usr/local/sbin/php-fpm
/usr/local/share/fpm/status.html
