# tinycorelinux/13.x/armv6/tcz/samba4.tcz.list
usr/local/etc/init.d/samba4
usr/local/share/samba4/COPYING
usr/local/share/samba4/files/smb.conf.default
usr/local/bin/mdfind
usr/local/bin/smbstatus
usr/local/bin/mvxattr
usr/local/bin/net
usr/local/bin/testparm
usr/local/bin/dbwrap_tool
usr/local/bin/ntlm_auth
usr/local/bin/sharesec
usr/local/bin/smbcquotas
usr/local/bin/smbcacls
usr/local/bin/nmblookup
usr/local/bin/smbget
usr/local/bin/pdbedit
usr/local/bin/smbpasswd
usr/local/bin/smbtree
usr/local/bin/smbcontrol
usr/local/bin/profiles
usr/local/bin/rpcclient
usr/local/bin/test_vfs_full_audit
usr/local/bin/test_nfs4_acls
usr/local/bin/vlp
usr/local/bin/rpc_open_tcp
usr/local/bin/timelimit
usr/local/bin/vfstest
usr/local/bin/pdbtest
usr/local/bin/test_mdsparser_es
usr/local/bin/smbconftort
usr/local/bin/smbtorture3
usr/local/bin/locktest2
usr/local/bin/smbspool_argv_wrapper
usr/local/bin/smbspool
usr/local/bin/smbclient
usr/local/bin/test_registry_regfio
usr/local/bin/test_tldap
usr/local/bin/pthreadpooltest_cmocka
usr/local/bin/pthreadpooltest
usr/local/bin/dumpmscat
usr/local/bin/dns_lookuptest
usr/local/bin/resolvconftest
usr/local/bin/test_gnutls
usr/local/bin/test_schannel
usr/local/bin/test_rc4_passwd_buffer
usr/local/bin/test_ntlm_check
usr/local/bin/test_ldap_message
usr/local/bin/async_connect_send_test
usr/local/bin/test_smb1cli_session
usr/local/bin/cifsdd
usr/local/bin/test_ndr_dns_nbt
usr/local/bin/test_ndr_string
usr/local/bin/locktest
usr/local/bin/masktest
usr/local/bin/gentest
usr/local/bin/smbtorture
usr/local/bin/oLschema2ldif
usr/local/bin/test_util_paths
usr/local/bin/test_byteorder
usr/local/bin/test_talloc_keep_secret
usr/local/bin/test_ms_fnmatch
usr/local/bin/test_rfc1738
usr/local/bin/regtree
usr/local/bin/regshell
usr/local/bin/regpatch
usr/local/bin/regdiff
usr/local/bin/wbinfo
usr/local/bin/stress-nss-libwbclient
usr/local/bin/nsstest
usr/local/bin/test_kerberos
usr/local/bin/test_encrypted_secrets_tdb
usr/local/bin/test_unique_object_sids
usr/local/bin/ndrdump
usr/local/bin/ldbrename
usr/local/bin/ldbedit
usr/local/bin/ldbmodify
usr/local/bin/ldbdel
usr/local/bin/ldbsearch
usr/local/bin/ldbadd
usr/local/bin/tdbtool
usr/local/bin/tdbbackup
usr/local/bin/tdbdump
usr/local/bin/tdbrestore
usr/local/bin/texpect
usr/local/bin/findsmb
usr/local/sbin/nmbd
usr/local/sbin/eventlogadm
usr/local/sbin/winbindd
usr/local/sbin/smbd
usr/local/sbin/samba-gpupdate
